# usr/lib/systemd/system/coreos-warn-invalid-mounts.service.d/10-mios-wsl2.conf
# /usr/lib/systemd/system/coreos-warn-invalid-mounts.service.d/10-mios-wsl2.conf
#
# coreos-warn-invalid-mounts.service prints a multi-line MOTD warning when /
# is not composefs-mounted and /sysroot is not read-only. Both conditions
# are STRUCTURALLY IMPOSSIBLE on WSL2: the WSL kernel mounts the distro
# rootfs as overlay via /init before systemd starts, with no composefs
# pathway and no read-only enforcement. On WSL the warning is purely noise
# and clutters every login.
#
# This drop-in skips the warning unit on WSL only. Bare-metal and VM
# deployments still see the warning if their root deployment drifts off
# the bootc/composefs path.
[Unit]
ConditionVirtualization=!wsl
