# usr/lib/systemd/system/cloud-final.service.d/10-mios-wsl2.conf
# Skip cloud-init's final stage in WSL.
#
# cloud-final.service runs `nc /run/cloud-init/share/final.sock` to
# signal completion to a hypervisor / cloud platform. WSL2 isn't a
# cloud platform; the socket is never created and the nc call fails
# with `nc: /run/cloud-init/share/final.sock: No such file or directory`
# at every boot, surfacing as the trailing error in the operator's
# `wsl -d MiOS` console output. Matches the existing skip pattern for
# cloud-init-local / cloud-init-network / cloud-config in the same
# directory neighbors.
[Unit]
ConditionVirtualization=!wsl
