# usr/lib/systemd/system/ublue-nvctk-cdi.service.d/10-mios.conf
# AI-hint: Configures the ublue-nvctk-cdi service to only activate if the NVIDIA device driver is present, ensuring the NVIDIA Container Toolkit CDI specification is only generated on systems with functional GPU hardware.
[Unit]
ConditionPathExists=/dev/nvidiactl