# usr/lib/systemd/system/mios-gpu-nvidia.service.d/10-cycle-fix.conf
# AI-hint: Overrides the mios-gpu-nvidia.service unit to ensure it runs during the early sysinit.target phase, bypassing standard dependency delays to ensure GPU drivers are initialized early in the boot sequence.
# AI-related: mios-gpu-nvidia, mios-gpu-nvidia.service, sysinit.target
[Unit]
DefaultDependencies=no
Requires=sysinit.target
After=sysinit.target
