# usr/share/containers/systemd/crowdsec-dashboard.container
[Unit]
Description='MiOS' CrowdSec Dashboard
After=network-online.target crowdsec.service
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips when CrowdSec isn't
# configured (no /etc/crowdsec/config.yaml = no LAPI to scrape).
ConditionPathExists=/etc/crowdsec/config.yaml

[Container]
Image=docker.io/crowdsecurity/crowdsec:latest
Network=mios.network
AutoUpdate=registry
User=mios-crowdsec
Group=mios-crowdsec

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
