# usr/share/containers/systemd/guacd.container
[Unit]
Description='MiOS' Apache Guacamole Daemon
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips inside nested
# containers (guacd needs raw VNC/RDP socket access).
ConditionVirtualization=!container

[Container]
# Image / Network placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml
# [image.sidecars].guacd and [network.quadlet].network.
Image=${MIOS_GUACD_IMAGE:-docker.io/guacamole/guacd:latest}
Network=${MIOS_QUADLET_NETWORK:-mios.network}
AutoUpdate=registry
User=mios-guacd
Group=mios-guacd

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
