# usr/share/containers/systemd/mios-crawl4ai.container
# /usr/share/containers/systemd/mios-crawl4ai.container
# 'MiOS' crawl4ai -- slim crawl engine (NO bundled Chromium).
#
# Built from usr/share/mios/crawl4ai/Containerfile (python:3.12 -- the host's
# 3.14 has no wheels for crawl4ai's deps -- + crawl4ai core + camoufox). The
# PRIMARY engine attaches to the EXISTING host Chrome over CDP; camoufox
# (stealth Firefox, baked in the image) is the fail-retry. The `crawl` broker
# verb (mios-crawl) POSTs to the loopback FastAPI service this container runs
# (mios-crawl4ai-service.py, mounted below).
[Unit]
Description='MiOS' crawl4ai slim engine (Chrome-CDP primary + camoufox fallback)
After=network-online.target mios-hermes-browser.service
Wants=network-online.target mios-hermes-browser.service

[Container]
# Slim local image (built by the .build step / podman build from the repo
# Containerfile). NO Chromium inside -- crawl4ai attaches to the host Chrome.
Image=${MIOS_CRAWL4AI_IMAGE:-localhost/mios-crawl4ai-slim:latest}
ContainerName=mios-crawl4ai
# Host networking: inside the container 127.0.0.1:9222 IS the host's ChromeDev
# CDP (the engine attaches there), and the service binds the host loopback
# :11235 directly -- same WSL host-net rationale as mios-ollama-cpu, and it
# avoids the podman-WSL PublishPort quirk.
Network=host
# Run as root INSIDE the container (isolated namespace): camoufox's patched
# Firefox is baked under root's cache at image-build time, so the service must
# read it as root. The container touches nothing host-writable except the RO
# service module below.
# The service module, mounted RO so it iterates without an image rebuild.
Volume=/usr/lib/mios/crawl4ai/mios-crawl4ai-service.py:/app/mios-crawl4ai-service.py:ro
Environment=MIOS_CRAWL_CDP_URL=${MIOS_CRAWL_CDP_URL:-http://127.0.0.1:9222}
Environment=MIOS_PORT_CRAWL4AI=${MIOS_PORT_CRAWL4AI:-11235}
Environment=MIOS_CRAWL_BIND=127.0.0.1
Environment=MIOS_CRAWL_CAMOUFOX=${MIOS_CRAWL_CAMOUFOX:-true}
Environment=MIOS_CRAWL_MIN_CHARS=${MIOS_CRAWL_MIN_CHARS:-200}
AutoUpdate=local
Label=org.opencontainers.image.title=mios-crawl4ai
Label=io.podman_desktop.openInBrowser=http://localhost:11235/healthz

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=120s

[Install]
WantedBy=multi-user.target default.target
