# usr/share/containers/systemd/mios-pxe-hub.container
[Unit]
Description='MiOS' PXE Boot Hub
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips on virtualized hosts
# without a routable LAN (PXE requires raw DHCP/TFTP on the host network).
ConditionVirtualization=!wsl
ConditionVirtualization=!container

[Container]
# Image / Network placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml
# [image.sidecars].pxe_hub and [network.quadlet].network.
Image=quay.io/poseidon/matchbox:latest
Network=mios.network
AutoUpdate=registry
User=mios-pxe-hub
Group=mios-pxe-hub

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
