# usr/share/containers/systemd/mios-guacamole.container
[Unit]
Description='MiOS' Apache Guacamole Web
After=network-online.target guacamole-postgres.service
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips when nested in another
# container (no working network namespacing for the bridge backend).
# guacamole-postgres must be reachable, hence After= on the DB unit.
ConditionVirtualization=!container

[Container]
# Image / Port placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml [image.sidecars].guacamole
# and [ports].guacamole_web.
Image=${MIOS_GUACAMOLE_IMAGE:-docker.io/guacamole/guacamole:latest}
Network=${MIOS_QUADLET_NETWORK:-mios.network}
AutoUpdate=registry
User=mios-guacamole
Group=mios-guacamole
# Tomcat inside the container listens 0.0.0.0:8080 by default. Map to
# host port guacamole_web (default 8090) because the localai port is
# already taken by mios-ai (LocalAI OpenAI-API endpoint, LAW 5). The
# /guacamole context path is preserved so links from the README and
# MOTD still resolve.
PublishPort=0.0.0.0:${MIOS_GUACAMOLE_PORT:-8090}:8080

# OCI labels for Podman Desktop discovery. Values MUST NOT contain
# spaces -- Quadlet's Label= parser word-splits on whitespace and breaks
# the generated podman-run command. Path includes /guacamole because
# the upstream image serves on that context root.
Label=org.opencontainers.image.title=mios-guacamole
Label=org.opencontainers.image.url=http://localhost:8090/guacamole/
Label=org.opencontainers.image.documentation=https://guacamole.apache.org/doc/
Label=io.podman_desktop.openInBrowser=http://localhost:8090/guacamole/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
