# usr/share/containers/systemd/mios-guacamole-postgres.container
# AI-hint: Systemd Quadlet configuration for the Guacamole PostgreSQL database container, defining the image, network, and lifecycle policies for the database backend in the MiOS environment.
# AI-related: mios-postgres, network-online.target, multi-user.target, default.target
[Unit]
Description='MiOS' Guacamole PostgreSQL
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips inside another
# container (postgres datadir on overlayfs-on-overlayfs is unreliable).
ConditionVirtualization=!container

[Container]
# Image / Network placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml
# [image.sidecars].postgres and [network.quadlet].network.
Image=docker.io/library/postgres:latest
ContainerName=mios-guacamole-postgres
Network=mios.network
AutoUpdate=registry
User=mios-postgres
Group=mios-postgres

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
