# usr/share/containers/systemd/mios-guacamole.container
[Unit]
Description='MiOS' Apache Guacamole Web
After=network-online.target guacamole-postgres.service
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips when nested in another
# container (no working network namespacing for the bridge backend).
# guacamole-postgres must be reachable, hence After= on the DB unit.
ConditionVirtualization=!container

[Container]
Image=docker.io/guacamole/guacamole:latest
Network=mios.network
AutoUpdate=registry
User=mios-guacamole
Group=mios-guacamole
# Tomcat inside the container listens 0.0.0.0:8080 by default. Map to
# host port 8090 because :8080 is already taken by mios-ai (LocalAI
# OpenAI-API endpoint -- Architectural Law 5). The /guacamole context
# path is preserved so links from the README and MOTD still resolve.
PublishPort=0.0.0.0:8090:8080

# OCI labels for Podman Desktop discovery -- "Browser desktop" entry
# point. Path includes /guacamole because the upstream image serves on
# that context root.
Label=org.opencontainers.image.title=MiOS Apache Guacamole (Browser desktop)
Label=org.opencontainers.image.description=Web client for RDP / VNC / SSH connections, served on port 8090 (mapped from container 8080 to avoid mios-ai conflict).
Label=org.opencontainers.image.url=http://localhost:8090/guacamole/
Label=org.opencontainers.image.documentation=https://guacamole.apache.org/doc/
Label=io.podman_desktop.openInBrowser=http://localhost:8090/guacamole/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
