# usr/share/containers/systemd/mios-llm-heavy.container
# AI-hint: GENERATED Quadlet container 'mios-llm-heavy' (WS-7 pods-as-SSOT). DO NOT EDIT -- regenerate via tools/generate-pod-quadlets.py from [containers.mios-llm-heavy] in mios.toml.
# /usr/share/containers/systemd/mios-llm-heavy.container

[Unit]
After=network-online.target mios-cdi-detect.service
ConditionPathExists=/usr/share/mios/vllm/model/config.json
Description='MiOS' vLLM heavy lane (OpenAI /v1, PagedAttention + APC; gated)
Wants=network-online.target

[Container]
AddDevice=nvidia.com/gpu=all
ContainerName=mios-llm-heavy
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=VLLM_NO_USAGE_STATS=1
Environment=DO_NOT_TRACK=1
Environment=VLLM_USE_V1=0
Environment=HF_HOME=/tmp/.cache
Environment=XDG_CACHE_HOME=/tmp/.cache
Environment=VLLM_CACHE_DIR=/tmp/.cache
Environment=FLASHINFER_WORKSPACE_DIR=/tmp/.cache
Environment=TRITON_CACHE_DIR=/tmp/.cache
Environment=VLLM_HOST_IP=127.0.0.1
Environment=GLOO_SOCKET_IFNAME=lo
Environment=NCCL_SOCKET_IFNAME=lo
Environment=HOME=/tmp
Exec=--model /models --served-model-name mios-heavy --host 0.0.0.0 --port ${MIOS_PORT_VLLM} --gpu-memory-utilization 0.45 --max-model-len 16384 --kv-cache-dtype auto --enable-prefix-caching --enforce-eager --enable-auto-tool-choice --tool-call-parser hermes
Group=815
Image=docker.io/vllm/vllm-openai:latest
Label=org.opencontainers.image.title=mios-llm-heavy
Label=org.opencontainers.image.url=http://localhost:11441/v1/models
Label=io.podman_desktop.openInBrowser=http://localhost:11441/v1/models
Pod=mios-ai.pod
User=815
Volume=/usr/share/mios/vllm/model:/models:ro,Z

[Install]
WantedBy=multi-user.target default.target

[Service]
Delegate=yes
EnvironmentFile=/etc/mios/install.env
Restart=on-failure
RestartSec=15
TimeoutStartSec=600s
