# usr/share/containers/systemd/mios-crowdsec-dashboard.container
# AI-hint: Systemd unit file defining the CrowdSec dashboard container deployment, managing the LAPI scraping service, volume mounts, and auto-update logic for the MiOS security monitoring stack.
# AI-related: crowdsec.service, network-online.target, multi-user.target, default.target
[Unit]
Description='MiOS' CrowdSec Dashboard
After=network-online.target crowdsec.service
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips when CrowdSec isn't
# configured (no /etc/crowdsec/config.yaml = no LAPI to scrape).
ConditionPathExists=/etc/crowdsec/config.yaml

[Container]
# Image / Network placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml
# [image.sidecars].crowdsec and [network.quadlet].network.
Image=docker.io/crowdsecurity/crowdsec:latest
ContainerName=mios-crowdsec-dashboard
Pod=mios-netinfra-sec.pod
AutoUpdate=registry
# Run as root to allow symlink creation and database management
# in the shared /var/lib/crowdsec volume.
User=0
Group=0

Volume=/var/lib/crowdsec:/var/lib/crowdsec:Z

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
