# usr/share/containers/systemd/mios-llm-worker@.container
# AI-hint: GENERATED Quadlet container 'mios-llm-worker@' (WS-7 pods-as-SSOT). DO NOT EDIT -- regenerate via tools/generate-pod-quadlets.py from [containers.mios-llm-worker@] in mios.toml.
# /usr/share/containers/systemd/mios-llm-worker@.container

[Unit]
After=network-online.target mios-cdi-detect.service
ConditionPathExists=/run/mios/swarm/%i.env
ConditionPathExists=/usr/share/mios/llamacpp/models/.ready
Description='MiOS' swarm worker %i (single-model llama-server, FOSS)
Wants=network-online.target

[Container]
AddDevice=nvidia.com/gpu=all
ContainerName=mios-llm-worker-%i
Entrypoint=/app/llama-server
Environment=HOME=/var/lib/mios/llamacpp
Environment=LD_LIBRARY_PATH=/usr/lib/wsl/lib:/usr/local/cuda/lib64
EnvironmentFile=-/run/mios/swarm/%i.env
Exec=--model /models/${WORKER_MODEL} --host 0.0.0.0 --port ${WORKER_PORT} --n-gpu-layers 99 --ctx-size 8192 --slot-save-path /var/lib/mios/llamacpp/slots/%i --parallel 1 --cont-batching --jinja --alias ${WORKER_MODEL}
Group=827
Image=ghcr.io/mostlygeek/llama-swap:cuda
Label=org.opencontainers.image.title=mios-llm-worker-%i
Pod=mios-ai.pod
User=827
Volume=/usr/share/mios/llamacpp/models:/models:ro,Z
Volume=/var/lib/mios/llamacpp/slots/%i:/var/lib/mios/llamacpp/slots/%i:Z

[Install]
WantedBy=multi-user.target default.target

[Service]
Delegate=yes
Restart=on-failure
RestartSec=10s
TimeoutStartSec=180s
