# usr/share/containers/systemd/mios-searxng.container
# AI-hint: Systemd container unit defining the SearXNG metasearch proxy; provides a privacy-preserving search tool for the LLM via port 8899, abstracting 250+ engines into a single local endpoint.
# AI-related: /etc/mios/searxng/, /etc/mios/searxng, /etc/mios/searxng/settings.yml, mios-searxng, mios-hermes, mios-services, network-online.target, multi-user.target, default.target, localhost:8899
# /etc/containers/systemd/mios-searxng.container
# 'MiOS' privacy-respecting metasearch (SearXNG upstream).
#
# Search counterpart to mios-hermes/mios-llm-light: a single localhost endpoint
# (for the operator and the LLM via tools) that
# aggregates results from ~250 upstream search engines without leaking
# any IP, user-agent, or session identifier to those engines. It is
# AGPL-3.0 and 100% local -- no external account, no API key, no
# telemetry. Pairs naturally with the agent surface as a `web_search`
# tool target without breaking Architectural Law 5 (no vendor cloud
# AI URL): SearXNG is a *search proxy*, not a model.
#
# Defaults policy: enabled by default. ConditionPathIsDirectory guards
# the config dir; if /etc/mios/searxng/ doesn't exist (incomplete
# bootstrap) the unit no-ops at pre-boot rather than crash-looping.

[Unit]
Description='MiOS' SearXNG metasearch (privacy-respecting search proxy)
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/etc/mios/searxng

[Container]
# Image / Port / UID placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml [image.sidecars].searxng,
# [ports].searxng, [services.searxng].uid/.gid. Multi-arch index manifest
# at the upstream tag covers amd64 / arm64 / armv7.
Image=docker.io/searxng/searxng:latest
ContainerName=mios-searxng
Pod=mios-webtools.pod
# The agent reaches SearXNG via SEARXNG_URL (host/localhost), not container
# DNS -- one network is enough (ai-net consolidated away 2026-05-23).
# SearXNG listens on 8080 inside the container.
# Map to host's loopback (it's consumed only by host hermes)
EnvironmentFile=/etc/mios/install.env

# Persistent state: settings.yml + uwsgi.ini live under /etc/searxng;
# the container expects them at /etc/searxng inside. We bind-mount the
# host's /etc/mios/searxng read-write so the operator (and SearXNG's
# own first-run secret_key generator) can persist edits across image
# rebuilds. /var/cache/mios/searxng holds the on-disk results cache so
# a restart doesn't repeat queries that the upstream engines have
# already answered.
Volume=/etc/mios/searxng:/etc/searxng:Z
Volume=/var/cache/mios/searxng:/var/cache/searxng:Z

# Identity. UID resolved from mios.toml [services.searxng].uid (default
# 818, pinned in usr/lib/sysusers.d/50-mios-services.conf). The upstream
# SearXNG image's `searxng` user inside the container is uid 977; we're
# overriding it so host-side bind-mounts are owned by mios-searxng and
# the container process can still read+write because the volumes are
# chowned <uid>:<gid> by usr/lib/tmpfiles.d/mios-searxng.conf at boot.
User=${MIOS_SEARXNG_UID:-818}
Group=${MIOS_SEARXNG_GID:-818}

# Public base URL for generated absolute links (RSS, opensearch.xml).
# Operators fronting SearXNG behind a reverse proxy override this in
# /etc/mios/searxng/settings.yml -> server.base_url.
Environment=SEARXNG_BASE_URL=http://localhost:${MIOS_PORT_SEARXNG:-8899}/
# SearXNG's entrypoint runs granian, which reads GRANIAN_HOST/PORT.
# Pin both to 8899 so host-network and bridge-network deployments
# both bind where MiOS expects. BIND_ADDRESS/SEARXNG_BIND_ADDRESS
# are kept for upstream Helm-chart parity.
Environment=SEARXNG_BIND_ADDRESS=0.0.0.0:${MIOS_PORT_SEARXNG:-8899}
Environment=BIND_ADDRESS=0.0.0.0:${MIOS_PORT_SEARXNG:-8899}
Environment=GRANIAN_HOST=0.0.0.0
Environment=GRANIAN_PORT=8899

# Concurrency (operator 2026-05-22: "SearXNG setup to handle the load").
# granian defaults to a single worker, which serialises concurrent incoming
# queries -- a bottleneck once the agent fans web_search out (per-query
# expansion x multi-agent council/DAG). SearXNG is a Flask WSGI app, so its
# request handlers run on granian BLOCKING threads. WORKERS x BLOCKING_THREADS
# is the concurrent-request ceiling; keep it >= MIOS_WEB_CONCURRENCY x
# MIOS_WEB_FANOUT (agent-pipe defaults 3 x 3 = 9, more under council/DAG
# fan-out). WORKERS=1 (1 x 8 = 8) sits BELOW even the single-agent peak, so a
# fan-out burst queued behind one worker and timed out ("Could not reach
# SearXNG ... timed out" in hermes-agent, 2026-06-04). 2 x 8 = 16 covers the
# peak with headroom. server.limiter stays false (settings.yml) so this
# self-querying burst is never rate-limited as bot traffic.
Environment=GRANIAN_WORKERS=2
Environment=GRANIAN_BLOCKING_THREADS=8

# Limit log volume; SearXNG is chatty per-query.
Environment=SEARXNG_LOG_LEVEL=WARNING

# OCI image labels for Podman Desktop. Values MUST NOT contain
# whitespace -- Quadlet's Label= parser word-splits on spaces and
# emits each word as a separate --label arg, breaking the run line.
Label=org.opencontainers.image.title=mios-searxng
Label=org.opencontainers.image.url=http://localhost:${MIOS_PORT_SEARXNG:-8899}/
Label=org.opencontainers.image.documentation=https://docs.searxng.org/
Label=io.podman_desktop.openInBrowser=http://localhost:${MIOS_PORT_SEARXNG:-8899}/

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
