# usr/share/containers/systemd/mios-cpu-node.container
# AI-hint: GENERATED Quadlet container 'mios-cpu-node' (WS-7 pods-as-SSOT). DO NOT EDIT -- regenerate via tools/generate-pod-quadlets.py from [containers.mios-cpu-node] in mios.toml.
# /usr/share/containers/systemd/mios-cpu-node.container

[Unit]
After=network-online.target
ConditionPathExists=/usr/share/mios/llamacpp/models/granite-4.1-8b.gguf
Description='MiOS' CPU inference node (always-on / gaming-immune last-resort lane)
Wants=network-online.target

[Container]
ContainerName=mios-cpu-node
Entrypoint=/app/llama-server
Environment=HOME=/var/lib/mios/llamacpp
Exec=--model /models/granite-4.1-8b.gguf --host 0.0.0.0 --port ${MIOS_PORT_CPU_NODE:-11451} --ctx-size 131072 --flash-attn on --cache-type-k q4_0 --cache-type-v q4_0 --parallel 1 --threads ${MIOS_CPU_NODE_THREADS:-14} --n-gpu-layers 0 --jinja --alias mios-cpu --slot-save-path ${MIOS_LLAMACPP_SLOT_DIR:-/var/lib/mios/llamacpp/slots}
Group=${MIOS_LLAMACPP_GID:-827}
Image=${MIOS_LLM_LIGHT_IMAGE:-ghcr.io/mostlygeek/llama-swap:cuda}
Label=org.opencontainers.image.title=mios-cpu-node
Pod=mios-ai.pod
User=${MIOS_LLAMACPP_UID:-827}
Volume=/usr/share/mios/llamacpp/models:/models:ro,Z
Volume=${MIOS_LLAMACPP_SLOT_DIR:-/var/lib/mios/llamacpp/slots}:${MIOS_LLAMACPP_SLOT_DIR:-/var/lib/mios/llamacpp/slots}:Z

[Install]
WantedBy=multi-user.target default.target

[Service]
Delegate=yes
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
