# usr/share/containers/systemd/mios-open-webui.container
# AI-hint: GENERATED Quadlet container 'mios-open-webui' (WS-7 pods-as-SSOT). DO NOT EDIT -- regenerate via tools/generate-pod-quadlets.py from [containers.mios-open-webui] in mios.toml.
# /usr/share/containers/systemd/mios-open-webui.container

[Unit]
After=network-online.target mios-gateway-agent.service mios-cdi-detect.service mios-webtools-firstboot.service
ConditionPathExists=/etc/mios/hermes/api.env
Description='MiOS' Open WebUI
Wants=network-online.target mios-gateway-agent.service mios-cdi-detect.service

[Container]
AutoUpdate=registry
ContainerName=mios-open-webui
Environment=OPENAI_API_BASE_URL=http://host.containers.internal:${MIOS_PORT_AGENT_PIPE}/v1
Environment=OLLAMA_BASE_URL=http://host.containers.internal:${MIOS_PORT_LLM_LIGHT}
Environment=WEBUI_NAME="MiOS AI"
Environment=WEBUI_SECRET_KEY=mios-stable-secret-change-me
Environment=ENABLE_SIGNUP=False
Environment=DEFAULT_USER_ROLE=admin
Environment=DEFAULT_MODELS=MiOS AI
Environment=PORT=${MIOS_PORT_OPEN_WEBUI}
Environment=ENABLE_WEB_SEARCH=True
Environment=ENABLE_IMAGE_GENERATION=True
Environment=ENABLE_CODE_INTERPRETER=True
Environment=ENABLE_API_KEYS=True
Environment=ENABLE_API_KEYS_ENDPOINT_RESTRICTIONS=False
Environment=ENABLE_FORWARD_USER_INFO_HEADERS=True
Environment=RAG_TOP_K=5
Environment=RAG_RELEVANCE_THRESHOLD=0.5
Environment=ENABLE_RAG_HYBRID_SEARCH=True
Environment=RAG_EMBEDDING_ENGINE=openai
Environment=RAG_EMBEDDING_MODEL=nomic-embed-text
Environment=RAG_OPENAI_API_BASE_URL=http://host.containers.internal:${MIOS_PORT_LLM_LIGHT}/v1
Environment=RAG_OPENAI_API_KEY=sk-noop
Environment=TASK_MODEL_EXTERNAL=hermes-agent
Environment=ENABLE_TITLE_GENERATION=True
Environment=ENABLE_TAGS_GENERATION=True
Environment=ENABLE_FOLLOW_UP_GENERATION=True
Environment=ENABLE_AUTOCOMPLETE_GENERATION=True
Environment="TITLE_GENERATION_PROMPT_TEMPLATE=Output a 3 to 6 word title for this chat, in the same language the chat used. Title only, no quotes, no preamble, no thinking. Chat:\n{{MESSAGES:END:2}}"
Environment="TAGS_GENERATION_PROMPT_TEMPLATE=Output 1 to 3 short topic tags for this chat (single words or hyphenated, in the chat language). JSON only: {\"tags\":[\"\",\"\"]}. Chat:\n{{MESSAGES:END:4}}"
Environment="FOLLOW_UP_GENERATION_PROMPT_TEMPLATE=Suggest 3 short follow up actions the user might want next, grounded in the CURRENT chat content and the operator's CURRENT system state. Each is a single phrase or short imperative the user could click, in the same language and tone the chat used. JSON only: {\"follow_ups\":[\"\",\"\",\"\"]}. Chat:\n{{MESSAGES:END:6}}"
Environment="AUTOCOMPLETE_GENERATION_PROMPT_TEMPLATE=Continue the user partial input naturally, in their language. Output ONLY the completion text, no quotes, no preamble. Input so far: {{PROMPT}}"
EnvironmentFile=/etc/mios/hermes/api.env
EnvironmentFile=/etc/mios/install.env
Exec=open-webui serve
Group=817
Image=localhost/mios-sys:latest
Label=org.opencontainers.image.title=mios-open-webui
Label=org.opencontainers.image.url=http://localhost:${MIOS_PORT_OPEN_WEBUI}/
Label=io.podman_desktop.openInBrowser=http://localhost:${MIOS_PORT_OPEN_WEBUI}/
Pod=mios-ai.pod
User=817
Volume=/var/lib/mios/open-webui:/app/backend/data:Z
Volume=/usr/share/mios:/usr/share/mios:ro
Volume=/etc/mios:/etc/mios:ro
Volume=/var/lib/mios/hermes-tail:/var/lib/mios/hermes-tail:ro
Volume=/var/lib/mios/delegation-prefilter:/var/lib/mios/delegation-prefilter:ro
Volume=/var/lib/mios/agent-nudger:/var/lib/mios/agent-nudger:ro
Volume=/var/lib/mios/log-watcher:/var/lib/mios/log-watcher:ro
Volume=/var/lib/mios/cron-director:/var/lib/mios/cron-director:ro
Volume=/usr/local/sbin:/usr/local/sbin:ro
Volume=/usr/libexec/mios:/usr/libexec/mios:ro
Volume=/var/lib/mios/agent-env:/var/lib/mios/agent-env:ro
Volume=/run/mios-launcher:/run/mios-launcher
Volume=/mnt/c:/mnt/c:ro
Volume=/mnt/m:/mnt/m:ro

[Install]
WantedBy=multi-user.target

[Service]
Delegate=yes
EnvironmentFile=/etc/mios/hermes/api.env
EnvironmentFile=/etc/mios/install.env
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
