# usr/share/containers/systemd/mios-ollama-igpu.container
[Unit]
Description='MiOS' Ollama (iGPU instance, micro-LLM lane via AMD ROCm)
After=network-online.target mios-cdi-detect.service
Wants=network-online.target
# Operator directive 2026-05-18: "iGPU model(s) isn't working" --
# diagnosis: ollama/ollama:latest is the CUDA-only image; even with
# the AMD CDI device (wsl2-amd.yaml) wired in via the
# /etc/containers/systemd/ollama.container.d/20-mios-igpu.conf
# drop-in, ollama-CUDA can't actually compute on the AMD iGPU. The
# fix is a SECOND ollama instance using ollama/ollama:rocm pinned to
# the AMD device. The original mios-ollama keeps the dGPU + CUDA;
# this instance gets the iGPU + ROCm. Sibling, not replacement.
#
# Operator directive 2026-05-17 (still in force): "iGPU's are ONLY
# micro-llms" + "wire to ollama!! JUST ONLY USES Micro-llms in MiOS
# stack". This Quadlet hosts only the small models (qwen3:1.7b,
# qwen3:0.6b, etc.) -- mios-daemon + agent-pipe router + (future)
# refine all point at :11435 for the micro-LLM lane. Big-model
# inference stays on the dGPU/CUDA mios-ollama at :11434.

[Container]
# ROCm variant -- only image that can compute on AMD GPUs. Pinned
# from mios.toml [image.sidecars].ollama_igpu (rendered at build).
Image=${MIOS_OLLAMA_IGPU_IMAGE:-docker.io/ollama/ollama:rocm}
ContainerName=mios-ollama-igpu
# Host networking -- same workaround as mios-surrealdb. A bridge
# (mios.network) + PublishPort hits a podman-WSL
# port-forwarder quirk where the host listener never binds even
# though `podman port` reports the mapping. Host-net gives the iGPU
# instance a reliable http://localhost:11435 for mios-daemon +
# agent-pipe router; no isolation loss that matters at loopback.
Network=host
AutoUpdate=registry

# AMD-only device pinning. Crucially does NOT add nvidia.com/gpu --
# this instance must stay off the dGPU so the dGPU/CUDA mios-ollama
# can keep using it for big models without contention. The amd.com/
# gpu CDI spec is the WSL2 hand-rolled wsl2-amd.yaml that
# mios-cdi-detect.service writes to /run/cdi/ at boot.
AddDevice=amd.com/gpu=all

# Identity: mios-ollama-igpu (uid 823) -- sibling sysuser to
# mios-ollama (uid 815). Separate UID so the iGPU instance's data
# dir ownership doesn't collide with the CUDA instance even on
# hosts that share /var/lib/ollama/models.
User=${MIOS_OLLAMA_IGPU_UID:-823}
Group=${MIOS_OLLAMA_IGPU_GID:-823}

# Distinct port -- mios-daemon's MIOS_DAEMON_ENDPOINT and the
# agent-pipe router's MIOS_AGENT_PIPE_ROUTER_ENDPOINT point at
# :11435 for the micro-LLM lane. The CUDA mios-ollama keeps :11434.
Environment=OLLAMA_HOST=0.0.0.0:${MIOS_PORT_OLLAMA_IGPU:-11435}
Environment=OLLAMA_ORIGINS=*

# Share the model store with mios-ollama -- the build-baked seed at
# /usr/share/ollama/models was hardlinked into /var/lib/ollama/models
# on first boot, and the iGPU instance reuses the same blobs. Saves
# disk + means a model pulled via either instance is immediately
# visible to the other.
Environment=OLLAMA_MODELS=/var/lib/ollama/models
Environment=HOME=/var/lib/ollama-igpu

# Micro-LLM lane parallelism. Smaller models, smaller per-request
# memory footprint -- can run more concurrent slots than the CUDA
# instance. mios-daemon's classify + refusal + cron loops all hit
# this lane; agent-pipe's router classifier hits it; future kanban-
# auto-dispatched parallel CPU agents will saturate it.
Environment=OLLAMA_NUM_PARALLEL=8
Environment=OLLAMA_MAX_LOADED_MODELS=6

# Volumes:
#   /var/lib/ollama          shared writable model store (read-write
#                             so the iGPU instance can ALSO pull
#                             models if asked; the small-model set
#                             rarely changes day-to-day)
#   /var/lib/ollama-igpu     iGPU-only HOME + ssh-key state
#   /usr/share/ollama        immutable seed mounted RO
Volume=/var/lib/ollama:/var/lib/ollama:Z
Volume=/var/lib/ollama-igpu:/var/lib/ollama-igpu:Z
Volume=/usr/share/ollama:/usr/share/ollama:ro,Z

# Host networking -- no PublishPort needed; the OLLAMA_HOST bind
# above (0.0.0.0:${MIOS_PORT_OLLAMA_IGPU}) goes directly to the
# host loopback + LAN interfaces.

Label=org.opencontainers.image.title=mios-ollama-igpu
Label=org.opencontainers.image.url=http://localhost:11435/
Label=io.podman_desktop.openInBrowser=http://localhost:11435/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
