# usr/share/containers/systemd/mios-ollama-cpu.container
[Unit]
Description='MiOS' Ollama (in-VM CPU light-lane: router/classifier + verb embeddings + micro-LLMs)
After=network-online.target mios-cdi-detect.service
Wants=network-online.target
# Honesty note 2026-05-24: this instance was historically named
# "mios-ollama-igpu" on the theory that ollama/ollama:rocm pinned to
# the AMD CDI device (wsl2-amd.yaml) would compute on the AMD iGPU.
# It never did -- ROCm cannot reach the iGPU inside WSL2 (the kernel
# doesn't expose /dev/kfd), so this instance has ALWAYS run on CPU.
# Renamed to mios-ollama-cpu to reflect that reality. The REAL AMD
# iGPU now runs natively on the Windows host
# (usr/share/mios/windows/mios-igpu-server.ps1) and is served to the
# swarm as mios-reasoner-cpu over Tailscale. The original mios-ollama
# keeps the dGPU + CUDA; this instance is the in-VM CPU light-lane.
# Sibling, not replacement.
#
# Operator directive 2026-05-17 (still in force): "iGPU's are ONLY
# micro-llms" + "wire to ollama!! JUST ONLY USES Micro-llms in MiOS
# stack". This Quadlet hosts only the small models (qwen3:1.7b,
# qwen3:0.6b, etc.) -- mios-daemon + agent-pipe router/classifier +
# verb embeddings + refine all point at :11435 for the CPU
# light-lane. Big-model inference stays on the dGPU/CUDA mios-ollama
# at :11434.

[Container]
# Standard (CPU/CUDA) image -- the rocm variant is pointless here
# since this instance runs on CPU only. Pinned from mios.toml
# [image.sidecars].ollama_cpu (rendered at build).
Image=${MIOS_OLLAMA_CPU_IMAGE:-docker.io/ollama/ollama:latest}
ContainerName=mios-ollama-cpu
# Host networking -- same workaround as mios-surrealdb. A bridge
# (mios.network) + PublishPort hits a podman-WSL
# port-forwarder quirk where the host listener never binds even
# though `podman port` reports the mapping. Host-net gives the CPU
# instance a reliable http://localhost:11435 for mios-daemon +
# agent-pipe router; no isolation loss that matters at loopback.
Network=host
AutoUpdate=registry

# No device pinning -- this instance runs on CPU only. It must stay
# off the dGPU so the dGPU/CUDA mios-ollama can keep using it for big
# models without contention, and there is no usable iGPU inside WSL2
# to attach (ROCm can't reach /dev/kfd here; the real AMD iGPU is the
# Windows-native mios-igpu-server.ps1 served as mios-reasoner-cpu).

# Identity: mios-ollama-cpu (uid 823) -- sibling sysuser to
# mios-ollama (uid 815). Separate UID so the CPU instance's data
# dir ownership doesn't collide with the CUDA instance even on
# hosts that share /var/lib/ollama/models.
User=${MIOS_OLLAMA_CPU_UID:-823}
Group=${MIOS_OLLAMA_CPU_GID:-823}

# Distinct port -- mios-daemon's MIOS_DAEMON_ENDPOINT and the
# agent-pipe router's MIOS_AGENT_PIPE_ROUTER_ENDPOINT point at
# :11435 for the CPU light-lane. The CUDA mios-ollama keeps :11434.
Environment=OLLAMA_HOST=0.0.0.0:${MIOS_PORT_OLLAMA_CPU:-11435}
Environment=OLLAMA_ORIGINS=*

# Share the model store with mios-ollama -- the build-baked seed at
# /usr/share/ollama/models was hardlinked into /var/lib/ollama/models
# on first boot, and the CPU instance reuses the same blobs. Saves
# disk + means a model pulled via either instance is immediately
# visible to the other.
Environment=OLLAMA_MODELS=/var/lib/ollama/models
Environment=HOME=/var/lib/ollama-cpu

# CPU light-lane parallelism. Smaller models, smaller per-request
# memory footprint -- can run more concurrent slots than the CUDA
# instance. mios-daemon's classify + refusal + cron loops all hit
# this lane; agent-pipe's router classifier hits it; future kanban-
# auto-dispatched parallel CPU agents will saturate it.
# RUNAWAY/OOM FIX (operator 2026-06-01): MAX_LOADED_MODELS was 6 with NO
# MemoryMax -- a wide research fan-out cold-loaded up to 6 distinct models
# here at once -> the kernel OOM killer fired repeatedly on this unit -> the
# whole VM wedged (3x terminate). Cap loaded models LOW (the light-lane only
# needs the micro-LLM resident at a time) and let parallel slots reuse it.
# SSOT: mios.toml [ai] knobs rendered to MIOS_OLLAMA_CPU_*.
Environment=OLLAMA_NUM_PARALLEL=${MIOS_OLLAMA_CPU_NUM_PARALLEL:-4}
Environment=OLLAMA_MAX_LOADED_MODELS=${MIOS_OLLAMA_CPU_MAX_LOADED:-2}

# Volumes:
#   /var/lib/ollama          shared writable model store (read-write
#                             so the CPU instance can ALSO pull
#                             models if asked; the small-model set
#                             rarely changes day-to-day)
#   /var/lib/ollama-cpu      CPU-instance HOME + ssh-key state
#   /usr/share/ollama        immutable seed mounted RO
Volume=/var/lib/ollama:/var/lib/ollama:Z
Volume=/var/lib/ollama-cpu:/var/lib/ollama-cpu:Z
Volume=/usr/share/ollama:/usr/share/ollama:ro,Z

# Host networking -- no PublishPort needed; the OLLAMA_HOST bind
# above (0.0.0.0:${MIOS_PORT_OLLAMA_CPU}) goes directly to the
# host loopback + LAN interfaces.

Label=org.opencontainers.image.title=mios-ollama-cpu
Label=org.opencontainers.image.url=http://localhost:11435/
Label=io.podman_desktop.openInBrowser=http://localhost:11435/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes
# OOM CONTAINMENT (operator 2026-06-01 wedge fix): cap this unit's memory so a
# cold-load storm OOM-kills WITHIN the cgroup instead of wedging the whole VM.
# MemoryHigh applies reclaim pressure first (soft); MemoryMax is the hard cap.
# Sized for the 40GiB-RAM dev VM: generous enough for normal micro-LLM load
# (a few small models resident) but a hard ceiling that contains a runaway
# herd to THIS cgroup instead of taking down the box. SSOT: mios.toml
# [ai].ollama_cpu_mem_{high,max} rendered to MIOS_OLLAMA_CPU_MEM_{HIGH,MAX}.
MemoryHigh=${MIOS_OLLAMA_CPU_MEM_HIGH:-12G}
MemoryMax=${MIOS_OLLAMA_CPU_MEM_MAX:-16G}

[Install]
WantedBy=multi-user.target default.target
