# usr/share/containers/systemd/ollama.container
[Unit]
Description='MiOS' Ollama LLM Server
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default. Ollama runs CPU-only when no
# GPU is present, so no virtualization-class condition is required.
# To force-disable, set [quadlets.enable] ollama = false in profile.toml.

[Container]
Image=docker.io/ollama/ollama:latest
Network=mios.network
AutoUpdate=registry
User=mios-ollama
Group=mios-ollama
# Bind 0.0.0.0 inside the container so the API is reachable from the
# host LAN, sibling containers, and loopback. PublishPort with explicit
# 0.0.0.0 on the host side makes the surface contract explicit.
Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=OLLAMA_ORIGINS=*
PublishPort=0.0.0.0:11434:11434

# OCI labels for Podman Desktop discovery. ollama is an alternate local
# LLM backend alongside mios-ai (LocalAI on :8080); MIOS_AI_ENDPOINT
# (Architectural Law 5) still points at LocalAI by default.
Label=org.opencontainers.image.title=MiOS Ollama LLM server
Label=org.opencontainers.image.description=Ollama API on port 11434 -- alternate local LLM runtime.
Label=org.opencontainers.image.url=http://localhost:11434/
Label=io.podman_desktop.openInBrowser=http://localhost:11434/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
