# usr/share/containers/systemd/ollama.container
[Unit]
Description='MiOS' Ollama LLM Server
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default. Ollama runs CPU-only when no
# GPU is present, so no virtualization-class condition is required.
# To force-disable, set [quadlets.enable] ollama = false in profile.toml.

[Container]
# Image / Port placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml [image.sidecars].ollama
# and [ports].ollama.
Image=${MIOS_OLLAMA_IMAGE:-docker.io/ollama/ollama:latest}
# ContainerName must be 'mios-ollama' -- ollama-firstboot.sh polls
# `podman exec mios-ollama ollama list` to know when the container is
# ready, and Podman Desktop's link discovery treats the name as the
# canonical handle. Quadlet's default would be `systemd-ollama`,
# which the firstboot script wouldn't find.
ContainerName=mios-ollama
# Dual-network: mios.network = historical core bridge; ai-net.network =
# AI-stack bridge so mios-hermes can reach mios-ollama:<port> by
# container name (deploy-aistack.sh's HERMES_BACKEND_BASE_URL target).
Network=${MIOS_QUADLET_NETWORK:-mios.network}
Network=ai-net.network
AutoUpdate=registry
# GPU passthrough. nvidia.com/gpu=all resolves to the CDI spec written
# by mios-cdi-detect.service (Before=ollama.service). On bare metal
# that's nvidia-ctk's generated /run/cdi/nvidia.yaml; on WSL2 it's the
# hand-rolled /run/cdi/wsl2-nvidia.yaml that maps /dev/dxg + rbinds
# /usr/lib/wsl. Operator-flagged 2026-05-11 (`load_tensors: CPU model
# buffer size = 16676 MiB` -- ollama was CPU-only because no CDI spec
# was attached). With this device line, ollama's runner.go reports
# `inference compute library=CUDA name=CUDA0 description=NVIDIA
# GeForce RTX 4090`. AMD/Intel hosts can swap the device class via
# /etc/containers/systemd/ollama.container.d/ drop-in overrides.
AddDevice=nvidia.com/gpu=all
# Numeric UID/GID -- the upstream ollama/ollama image has no
# `mios-ollama` user in its /etc/passwd, so a name-based User= lookup
# fails with "unable to find user mios-ollama: no matching entries in
# passwd file" (status 126). UID/GID resolved from mios.toml
# [services.ollama].uid/.gid (default 815) and pinned in
# usr/lib/sysusers.d/50-mios-services.conf so the host's mios-ollama
# user has the same UID and the /var/lib/ollama volume ownership
# stays consistent.
User=${MIOS_OLLAMA_UID:-815}
Group=${MIOS_OLLAMA_GID:-815}
# Bind 0.0.0.0 so the WSL VM's eth0 (NAT IP) accepts inbound from the
# Windows-host netsh portproxy that build-mios.ps1 sets up
# (0.0.0.0:11434 -> <WSL-VM-IP>:11434). Pre-2026-05-13 history: was
# 127.0.0.1 because the original design relied on WSL2 native
# localhostForwarding (which only catches AF_INET binds, not the
# AF_INET6 dual-stack `*:N` that ollama produces from 0.0.0.0).
# But the native relay turned out to silently fail for many ports
# under NAT mode -- explicit portproxy is now the canonical Windows-
# host -> WSL service path, and that requires the service to bind
# the WSL VM's eth0, not just container loopback. 0.0.0.0 also
# restores LAN-side reach (phone / tablet using local LLM via
# Open WebUI). Operator-confirmed 2026-05-13: with HOST=127.0.0.1,
# Windows-host localhost:11434 returned 000 even with portproxy +
# Hyper-V firewall opened.
Environment=OLLAMA_HOST=0.0.0.0:${MIOS_PORT_OLLAMA:-11434}
Environment=OLLAMA_ORIGINS=*
# Point ollama at the writable runtime store. mios-ollama-firstboot.
# service hardlink-copies the build-baked seed (/usr/share/ollama/
# models -- immutable composefs surface) into this path on first boot.
# Subsequent 'ollama pull' commands write here too, so this is the
# single source of truth ollama reads.
Environment=OLLAMA_MODELS=/var/lib/ollama/models

# Volumes (Architectural Law 2 -- /var paths declared via tmpfiles.d):
#   /var/lib/ollama   writable runtime model store
#   /usr/share/ollama immutable seed mounted RO so the firstboot service
#                     and ollama itself can both reach it
Volume=/var/lib/ollama:/var/lib/ollama:Z
Volume=/usr/share/ollama:/usr/share/ollama:ro,Z

PublishPort=0.0.0.0:${MIOS_PORT_OLLAMA:-11434}:11434

# OCI labels for Podman Desktop discovery. Values MUST NOT contain
# spaces -- Quadlet's Label= parser word-splits on whitespace and emits
# each word as a separate --label arg, breaking the generated
# podman-run command line. ollama is THE local model + embedding
# backend; Hermes-Agent (the LIVE MiOS agent at root) fronts it on
# :8642 as the canonical MIOS_AI_ENDPOINT per Architectural Law 5.
Label=org.opencontainers.image.title=mios-ollama
Label=org.opencontainers.image.url=http://localhost:11434/
Label=io.podman_desktop.openInBrowser=http://localhost:11434/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
