# usr/share/containers/systemd/ollama.container
[Unit]
Description='MiOS' Ollama LLM Server
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default. Ollama runs CPU-only when no
# GPU is present, so no virtualization-class condition is required.
# To force-disable, set [quadlets.enable] ollama = false in profile.toml.

[Container]
Image=docker.io/ollama/ollama:latest
Network=mios.network
AutoUpdate=registry
User=mios-ollama
Group=mios-ollama

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
