# usr/share/containers/systemd/ollama.container
[Unit]
Description='MiOS' Ollama LLM Server
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default. Ollama runs CPU-only when no
# GPU is present, so no virtualization-class condition is required.
# To force-disable, set [quadlets.enable] ollama = false in profile.toml.

[Container]
Image=docker.io/ollama/ollama:latest@sha256:5600a652d1081050f398152127c584222546354491f27fe47ccbc6351bc870bd
# ContainerName must be 'mios-ollama' -- ollama-firstboot.sh polls
# `podman exec mios-ollama ollama list` to know when the container is
# ready, and Podman Desktop's link discovery treats the name as the
# canonical handle. Quadlet's default would be `systemd-ollama`,
# which the firstboot script wouldn't find.
ContainerName=mios-ollama
# Dual-network: mios.network = historical core bridge; ai-net.network =
# AI-stack bridge so mios-hermes can reach mios-ollama:11434 by
# container name (deploy-aistack.sh's HERMES_BACKEND_BASE_URL target).
Network=mios.network
Network=ai-net.network
AutoUpdate=registry
User=mios-ollama
Group=mios-ollama
# Bind 0.0.0.0 inside the container so the API is reachable from the
# host LAN, sibling containers, and loopback. PublishPort with explicit
# 0.0.0.0 on the host side makes the surface contract explicit.
Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=OLLAMA_ORIGINS=*
# Point ollama at the writable runtime store. mios-ollama-firstboot.
# service hardlink-copies the build-baked seed (/usr/share/ollama/
# models -- immutable composefs surface) into this path on first boot.
# Subsequent 'ollama pull' commands write here too, so this is the
# single source of truth ollama reads.
Environment=OLLAMA_MODELS=/var/lib/ollama/models

# Volumes (Architectural Law 2 -- /var paths declared via tmpfiles.d):
#   /var/lib/ollama   writable runtime model store
#   /usr/share/ollama immutable seed mounted RO so the firstboot service
#                     and ollama itself can both reach it
Volume=/var/lib/ollama:/var/lib/ollama:Z
Volume=/usr/share/ollama:/usr/share/ollama:ro,Z

PublishPort=0.0.0.0:11434:11434

# OCI labels for Podman Desktop discovery. Values MUST NOT contain
# spaces -- Quadlet's Label= parser word-splits on whitespace and emits
# each word as a separate --label arg, breaking the generated
# podman-run command line. ollama is an alternate local LLM backend
# alongside mios-ai (LocalAI on :8080); MIOS_AI_ENDPOINT
# (Architectural Law 5) still points at LocalAI by default.
Label=org.opencontainers.image.title=mios-ollama
Label=org.opencontainers.image.url=http://localhost:11434/
Label=io.podman_desktop.openInBrowser=http://localhost:11434/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
