# usr/share/containers/systemd/qdrant.container
# /etc/containers/systemd/qdrant.container

[Unit]
Description=Qdrant vector database for MiOS KB
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=qdrant
Image=docker.io/qdrant/qdrant:latest
AutoUpdate=registry
Volume=/var/lib/mios/qdrant:/qdrant/storage:Z
# Bind 0.0.0.0 so WSL2 NAT-mode localhostForwarding reaches Windows host.
PublishPort=6333:6333
PublishPort=6334:6334
Environment=QDRANT__SERVICE__API_KEY=${QDRANT_API_KEY:-}

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
