# usr/share/containers/systemd/qdrant.container
# /etc/containers/systemd/qdrant.container

[Unit]
Description=Qdrant vector database for MiOS KB
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=qdrant
Image=docker.io/qdrant/qdrant:latest
AutoUpdate=registry
Volume=/var/lib/mios/qdrant:/qdrant/storage:Z
# Bind 0.0.0.0 so WSL2 NAT-mode localhostForwarding reaches Windows host.
PublishPort=6333:6333
PublishPort=6334:6334
Environment=QDRANT__SERVICE__API_KEY=${QDRANT_API_KEY:-}
# LAW 6 (UNPRIVILEGED-QUADLETS) compliance per AGENTS.md §5: every
# Quadlet must declare User=, Group=, Delegate=. qdrant needs root for
# the /qdrant/storage bind volume (upstream image runs as root and
# chown-s storage at startup). Documented Law 6 exception, same pattern
# as mios-hermes-workspace. Operator-confirmed 2026-05-13: build fails
# at 99-postcheck.sh "missing User= directive" without this.
User=0
Group=0
Delegate=yes

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
