# usr/lib/systemd/system/cockpit-wsinstance-https-factory.socket.d/10-mios-container.conf
# /usr/lib/systemd/system/cockpit-wsinstance-https-factory.socket.d/10-mios-container.conf
#
# Pair drop-in to cockpit-wsinstance-http.socket.d sibling. Same fix
# pattern: strip namespace requirements so the socket can bind inside
# a container / podman-machine WSL distro.

[Unit]
ConditionVirtualization=container

[Socket]
DynamicUser=no
PrivateTmp=no
PrivateMounts=no
PrivateDevices=no
PrivateUsers=no
ProtectHome=no
ProtectSystem=no
RestrictNamespaces=no
SocketUser=cockpit-ws
SocketGroup=cockpit-ws
