# usr/lib/systemd/system/cockpit-wsinstance-https-factory.socket.d/10-mios-container.conf
# /usr/lib/systemd/system/cockpit-wsinstance-https-factory.socket.d/10-mios-container.conf
#
# Pair drop-in to cockpit-wsinstance-http.socket.d sibling. Same fix
# pattern: strip namespace requirements so the socket can bind inside
# a container / podman-machine WSL distro.
#
# KEEP the stock SocketUser=root + SocketGroup=cockpit-wsinstance-socket
# (do NOT override to cockpit-ws:cockpit-ws). cockpit.service runs as
# User=cockpit-systemd-service with Group=cockpit-wsinstance-socket, so
# the cockpit-tls helper (which is what connects to this socket) has
# group-rw access only when SocketGroup matches its supplementary group.
# Operator-flagged 2026-05-10: overriding to cockpit-ws:cockpit-ws caused
#     cockpit-tls: connect(https-factory.sock) failed: Permission denied
# because cockpit-systemd-service isn't in group cockpit-ws.

[Unit]
ConditionVirtualization=container

[Socket]
DynamicUser=no
PrivateTmp=no
PrivateMounts=no
PrivateDevices=no
PrivateUsers=no
ProtectHome=no
ProtectSystem=no
RestrictNamespaces=no
