# usr/lib/systemd/system/mios-daemon.service.d/20-cpu-limit.conf
# /usr/lib/systemd/system/mios-daemon.service.d/20-cpu-limit.conf
# Operator 2026-05-18: keep the always-on micro-LLM daemon to "minimal
# background impact" -- cap TOTAL CPU at 50% (half a core's worth of
# CPU-seconds per real second). The daemon ALSO pins itself to 2
# randomly-chosen cores at startup (os.sched_setaffinity in mios-daemon
# main()), so it is confined to 2 cores AND throttled to 50% overall.
# CPUWeight + Nice make it yield to interactive / agent work.
[Service]
CPUQuota=50%
CPUWeight=20
Nice=10
