# usr/lib/systemd/system/mios-daemon.service.d/20-cpu-limit.conf
# /usr/lib/systemd/system/mios-daemon.service.d/20-cpu-limit.conf
# Operator 2026-05-24: the daemon-agent runs ALWAYS-ON background work
# (global journals/logs/system management) at ~2c/4t, but must BURST to 8c/16t
# when it participates in a FOREGROUND AI pipeline -- concurrently with the
# GPU (4090) + iGPU compute. So allow up to 8 cores (CPUQuota=800%) and pin to
# 8 cores (os.sched_setaffinity in mios-daemon main()). CPUWeight=20 + Nice=10
# keep it LIGHT at baseline and make it YIELD those cores to the GPU agents /
# interactive turns when contended -- the ~2c baseline is emergent from the
# light background loop + low weight; 8c is the foreground ceiling.
# (Was CPUQuota=50% + a 2-core pin -- the old "minimal background impact" era.)
[Service]
CPUQuota=800%
CPUWeight=20
Nice=10
