# usr/lib/systemd/system/polkit.service.d/10-mios-container.conf
# AI-hint: Configures polkit.service restart limits and retry intervals to ensure service stability and availability within MiOS containerized environments.
# AI-related: polkit.service
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Restart=on-failure
RestartSec=30
