# usr/lib/systemd/system/polkit.service.d/10-mios-container.conf
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Restart=on-failure
RestartSec=30
