# usr/lib/systemd/system/gnome-remote-desktop.service.d/10-network-wait.conf
# AI-hint: Ensures the gnome-remote-desktop service waits for a stable network-online.target before starting to prevent connection failures during boot.
# AI-related: network-online.target
[Unit]
After=network-online.target
Wants=network-online.target
