# usr/share/mios/templates/conformance-grandfathered.list
usr/lib/mios/agent-pipe/mios_pipe/__init__.py
usr/lib/mios/agent-pipe/mios_pipe/redact.py
usr/lib/mios/agent-pipe/test_mios_a2a.py
usr/lib/mios/agent-pipe/test_mios_a2a_client.py
usr/lib/mios/agent-pipe/test_mios_a2a_loopback.py
usr/lib/mios/agent-pipe/test_mios_a2a_passport.py
usr/lib/mios/agent-pipe/test_mios_a2a_principal.py
usr/lib/mios/agent-pipe/test_mios_aci.py
usr/lib/mios/agent-pipe/test_mios_admission.py
usr/lib/mios/agent-pipe/test_mios_agent_call.py
usr/lib/mios/agent-pipe/test_mios_agentreg.py
usr/lib/mios/agent-pipe/test_mios_ai_manifest.py
usr/lib/mios/agent-pipe/test_mios_antifab.py
usr/lib/mios/agent-pipe/test_mios_applet_webresearch.py
usr/lib/mios/agent-pipe/test_mios_approutes.py
usr/lib/mios/agent-pipe/test_mios_arbiter.py
usr/lib/mios/agent-pipe/test_mios_audit.py
usr/lib/mios/agent-pipe/test_mios_backfill.py
usr/lib/mios/agent-pipe/test_mios_batch.py
usr/lib/mios/agent-pipe/test_mios_bench.py
usr/lib/mios/agent-pipe/test_mios_bench_harness.py
usr/lib/mios/agent-pipe/test_mios_blades.py
usr/lib/mios/agent-pipe/test_mios_budget.py
usr/lib/mios/agent-pipe/test_mios_build_catalog.py
usr/lib/mios/agent-pipe/test_mios_capreg.py
usr/lib/mios/agent-pipe/test_mios_chat.py
usr/lib/mios/agent-pipe/test_mios_classify.py
usr/lib/mios/agent-pipe/test_mios_clusterhealth.py
usr/lib/mios/agent-pipe/test_mios_codemode.py
usr/lib/mios/agent-pipe/test_mios_cold_evict.py
usr/lib/mios/agent-pipe/test_mios_compact.py
usr/lib/mios/agent-pipe/test_mios_compound.py
usr/lib/mios/agent-pipe/test_mios_config.py
usr/lib/mios/agent-pipe/test_mios_config_audit.py
usr/lib/mios/agent-pipe/test_mios_cost.py
usr/lib/mios/agent-pipe/test_mios_council_diversity.py
usr/lib/mios/agent-pipe/test_mios_crl.py
usr/lib/mios/agent-pipe/test_mios_ctxpack.py
usr/lib/mios/agent-pipe/test_mios_cua.py
usr/lib/mios/agent-pipe/test_mios_cua_hierarchy.py
usr/lib/mios/agent-pipe/test_mios_daemon.py
usr/lib/mios/agent-pipe/test_mios_daemons.py
usr/lib/mios/agent-pipe/test_mios_dag_exec.py
usr/lib/mios/agent-pipe/test_mios_db_config.py
usr/lib/mios/agent-pipe/test_mios_dci.py
usr/lib/mios/agent-pipe/test_mios_dispatch.py
usr/lib/mios/agent-pipe/test_mios_dispatcher.py
usr/lib/mios/agent-pipe/test_mios_dual_ledger.py
usr/lib/mios/agent-pipe/test_mios_egress.py
usr/lib/mios/agent-pipe/test_mios_embed_backfill.py
usr/lib/mios/agent-pipe/test_mios_endpoints.py
usr/lib/mios/agent-pipe/test_mios_evict.py
usr/lib/mios/agent-pipe/test_mios_fanout.py
usr/lib/mios/agent-pipe/test_mios_firewall.py
usr/lib/mios/agent-pipe/test_mios_gateway_agent.py
usr/lib/mios/agent-pipe/test_mios_gateway_queue.py
usr/lib/mios/agent-pipe/test_mios_gossip.py
usr/lib/mios/agent-pipe/test_mios_grounding.py
usr/lib/mios/agent-pipe/test_mios_hitl.py
usr/lib/mios/agent-pipe/test_mios_hitlflow.py
usr/lib/mios/agent-pipe/test_mios_hopbudget.py
usr/lib/mios/agent-pipe/test_mios_http_caps.py
usr/lib/mios/agent-pipe/test_mios_interop.py
usr/lib/mios/agent-pipe/test_mios_jsonsalvage.py
usr/lib/mios/agent-pipe/test_mios_k3s.py
usr/lib/mios/agent-pipe/test_mios_kernel.py
usr/lib/mios/agent-pipe/test_mios_knowledge.py
usr/lib/mios/agent-pipe/test_mios_kvfork.py
usr/lib/mios/agent-pipe/test_mios_kvgc.py
usr/lib/mios/agent-pipe/test_mios_lanes.py
usr/lib/mios/agent-pipe/test_mios_lanes_resolver.py
usr/lib/mios/agent-pipe/test_mios_launch.py
usr/lib/mios/agent-pipe/test_mios_letta.py
usr/lib/mios/agent-pipe/test_mios_list_dir.py
usr/lib/mios/agent-pipe/test_mios_manifest.py
usr/lib/mios/agent-pipe/test_mios_mcp.py
usr/lib/mios/agent-pipe/test_mios_mcp_pool.py
usr/lib/mios/agent-pipe/test_mios_mcp_sandbox.py
usr/lib/mios/agent-pipe/test_mios_memguard.py
usr/lib/mios/agent-pipe/test_mios_memory.py
usr/lib/mios/agent-pipe/test_mios_mtls.py
usr/lib/mios/agent-pipe/test_mios_native_loop.py
usr/lib/mios/agent-pipe/test_mios_oscontrol.py
usr/lib/mios/agent-pipe/test_mios_owui.py
usr/lib/mios/agent-pipe/test_mios_pdp.py
usr/lib/mios/agent-pipe/test_mios_pg.py
usr/lib/mios/agent-pipe/test_mios_planner.py
usr/lib/mios/agent-pipe/test_mios_policy.py
usr/lib/mios/agent-pipe/test_mios_portal.py
usr/lib/mios/agent-pipe/test_mios_preempt.py
usr/lib/mios/agent-pipe/test_mios_promptfmt.py
usr/lib/mios/agent-pipe/test_mios_promptver.py
usr/lib/mios/agent-pipe/test_mios_provider_translate.py
usr/lib/mios/agent-pipe/test_mios_quarantine.py
usr/lib/mios/agent-pipe/test_mios_quota.py
usr/lib/mios/agent-pipe/test_mios_react_reflexion.py
usr/lib/mios/agent-pipe/test_mios_record_replay.py
usr/lib/mios/agent-pipe/test_mios_redact.py
usr/lib/mios/agent-pipe/test_mios_refine.py
usr/lib/mios/agent-pipe/test_mios_reflect.py
usr/lib/mios/agent-pipe/test_mios_registry.py
usr/lib/mios/agent-pipe/test_mios_reputation.py
usr/lib/mios/agent-pipe/test_mios_router.py
usr/lib/mios/agent-pipe/test_mios_routing.py
usr/lib/mios/agent-pipe/test_mios_ruleof2.py
usr/lib/mios/agent-pipe/test_mios_sandbox.py
usr/lib/mios/agent-pipe/test_mios_sched.py
usr/lib/mios/agent-pipe/test_mios_scratchpad.py
usr/lib/mios/agent-pipe/test_mios_secondary_loop.py
usr/lib/mios/agent-pipe/test_mios_secset.py
usr/lib/mios/agent-pipe/test_mios_selfimprove.py
usr/lib/mios/agent-pipe/test_mios_selfimprove_act.py
usr/lib/mios/agent-pipe/test_mios_skills.py
usr/lib/mios/agent-pipe/test_mios_slo.py
usr/lib/mios/agent-pipe/test_mios_smartroute.py
usr/lib/mios/agent-pipe/test_mios_sse.py
usr/lib/mios/agent-pipe/test_mios_stress.py
usr/lib/mios/agent-pipe/test_mios_surface.py
usr/lib/mios/agent-pipe/test_mios_swarm.py
usr/lib/mios/agent-pipe/test_mios_tiered_memory.py
usr/lib/mios/agent-pipe/test_mios_tokenize.py
usr/lib/mios/agent-pipe/test_mios_toolconflict.py
usr/lib/mios/agent-pipe/test_mios_toolexec.py
usr/lib/mios/agent-pipe/test_mios_toolsearch.py
usr/lib/mios/agent-pipe/test_mios_trace.py
usr/lib/mios/agent-pipe/test_mios_turn.py
usr/lib/mios/agent-pipe/test_mios_user_config.py
usr/lib/mios/agent-pipe/test_mios_vector.py
usr/lib/mios/agent-pipe/test_mios_verbcatalog.py
usr/lib/mios/agent-pipe/test_mios_verity.py
usr/lib/mios/agent-pipe/test_mios_vision.py
usr/lib/mios/agent-pipe/test_mios_web_research.py
usr/lib/mios/agent-pipe/test_mios_worker_tools.py
usr/lib/systemd/system/ceph-bootstrap.service
usr/lib/systemd/system/hermes-dashboard.service
usr/lib/systemd/system/hermes-worker-firstboot.service
usr/lib/systemd/system/hermes-worker.service
usr/lib/systemd/system/k3s.service
usr/lib/systemd/system/mios-additionalimagestores-perms.service
usr/lib/systemd/system/mios-adguard-firstboot.service
usr/lib/systemd/system/mios-agent-pipe.service
usr/lib/systemd/system/mios-agents.service
usr/lib/systemd/system/mios-ai-firstboot.service
usr/lib/systemd/system/mios-aios-refresh.service
usr/lib/systemd/system/mios-boot-diag.service
usr/lib/systemd/system/mios-bootc-switch.service
usr/lib/systemd/system/mios-cdi-detect.service
usr/lib/systemd/system/mios-ceph-bootstrap.service
usr/lib/systemd/system/mios-copy-build-log.service
usr/lib/systemd/system/mios-cpu-isolate.service
usr/lib/systemd/system/mios-cron-director.service
usr/lib/systemd/system/mios-daemon.service
usr/lib/systemd/system/mios-dashboard-issue.service
usr/lib/systemd/system/mios-delegation-prefilter.service
usr/lib/systemd/system/mios-embed-backfill.service
usr/lib/systemd/system/mios-finetune-serve.service
usr/lib/systemd/system/mios-firewall-ports.service
usr/lib/systemd/system/mios-flatpak-init.service
usr/lib/systemd/system/mios-flatpak-install.service
usr/lib/systemd/system/mios-forge-firstboot.service
usr/lib/systemd/system/mios-forgejo-runner-firstboot.service
usr/lib/systemd/system/mios-freeipa-enroll.service
usr/lib/systemd/system/mios-gateway-agent.service
usr/lib/systemd/system/mios-git-root-init.service
usr/lib/systemd/system/mios-gpu-amd.service
usr/lib/systemd/system/mios-gpu-detect.service
usr/lib/systemd/system/mios-gpu-intel.service
usr/lib/systemd/system/mios-gpu-nvidia.service
usr/lib/systemd/system/mios-gpu-pv-detect.service
usr/lib/systemd/system/mios-gpu-status.service
usr/lib/systemd/system/mios-grd-setup.service
usr/lib/systemd/system/mios-ha-bootstrap.service
usr/lib/systemd/system/mios-hermes-browser-worker.service
usr/lib/systemd/system/mios-hermes-browser.service
usr/lib/systemd/system/mios-hermes-firstboot.service
usr/lib/systemd/system/mios-hermes-tail.service
usr/lib/systemd/system/mios-hyperv-enhanced.service
usr/lib/systemd/system/mios-k3s-init.service
usr/lib/systemd/system/mios-kvmfr-load.service
usr/lib/systemd/system/mios-libexec-perms.service
usr/lib/systemd/system/mios-libvirtd-setup.service
usr/lib/systemd/system/mios-mcp.service
usr/lib/systemd/system/mios-opencode-gateway.service
usr/lib/systemd/system/mios-passport-provision.service
usr/lib/systemd/system/mios-pgvector-backup.service
usr/lib/systemd/system/mios-podman-gc.service
usr/lib/systemd/system/mios-podman-ps.service
usr/lib/systemd/system/mios-policy-arbiter.service
usr/lib/systemd/system/mios-role.service
usr/lib/systemd/system/mios-selinux-init.service
usr/lib/systemd/system/mios-skills-miner.service
usr/lib/systemd/system/mios-sriov-init.service
usr/lib/systemd/system/mios-suggestion-refresh.service
usr/lib/systemd/system/mios-swarm-pack-firstboot.service
usr/lib/systemd/system/mios-sync-theme.service
usr/lib/systemd/system/mios-sys-env-refresh.service
usr/lib/systemd/system/mios-ttyd-bash.service
usr/lib/systemd/system/mios-ttyd-expose.service
usr/lib/systemd/system/mios-ttyd-powershell.service
usr/lib/systemd/system/mios-verify-root.service
usr/lib/systemd/system/mios-verify.service
usr/lib/systemd/system/mios-waydroid-init.service
usr/lib/systemd/system/mios-webtools-firstboot.service
usr/lib/systemd/system/mios-wsl-early.service
usr/lib/systemd/system/mios-wsl-firstboot.service
usr/lib/systemd/system/mios-wsl-flatpak-export-sync.service
usr/lib/systemd/system/mios-wsl-init.service
usr/lib/systemd/system/mios-wsl-interop-priority.service
usr/lib/systemd/system/mios-wsl-runtime-dir.service
usr/lib/systemd/system/mios-wslg-permissions-fix.service
usr/libexec/mios/mios-a2a-delegate
usr/libexec/mios/mios-a2a-discover
usr/libexec/mios/mios-a2a-mdns
usr/libexec/mios/mios-a2a-test
usr/libexec/mios/mios-adguard-firstboot
usr/libexec/mios/mios-agents-firstboot.sh
usr/libexec/mios/mios-ai-capabilities-gen
usr/libexec/mios/mios-ai-clear
usr/libexec/mios/mios-ai-firstboot
usr/libexec/mios/mios-ai-hint-coverage
usr/libexec/mios/mios-ai-manifest-gen
usr/libexec/mios/mios-ai-reset
usr/libexec/mios/mios-ai-tag
usr/libexec/mios/mios-app-default
usr/libexec/mios/mios-app-search
usr/libexec/mios/mios-app-type
usr/libexec/mios/mios-apps
usr/libexec/mios/mios-as-operator
usr/libexec/mios/mios-autocenter
usr/libexec/mios/mios-bake-group
usr/libexec/mios/mios-bench
usr/libexec/mios/mios-build-driver
usr/libexec/mios/mios-build-status
usr/libexec/mios/mios-build-tail
usr/libexec/mios/mios-cache-clear
usr/libexec/mios/mios-cdi-detect
usr/libexec/mios/mios-cdp-fetch
usr/libexec/mios/mios-ceph-configure
usr/libexec/mios/mios-cephfs-provision
usr/libexec/mios/mios-chain-verify
usr/libexec/mios/mios-codebase-index
usr/libexec/mios/mios-codemode-api.py
usr/libexec/mios/mios-coderun
usr/libexec/mios/mios-coderun-broker
usr/libexec/mios/mios-coderun-codemode
usr/libexec/mios/mios-coderun-session
usr/libexec/mios/mios-compact
usr/libexec/mios/mios-computer-use
usr/libexec/mios/mios-computer-use-server
usr/libexec/mios/mios-configurator-launch
usr/libexec/mios/mios-crawl
usr/libexec/mios/mios-cron-director
usr/libexec/mios/mios-cron-schedule
usr/libexec/mios/mios-cursor-apply
usr/libexec/mios/mios-cursor-ensure
usr/libexec/mios/mios-daemon
usr/libexec/mios/mios-dashboard-render-issue.sh
usr/libexec/mios/mios-dashboard.sh
usr/libexec/mios/mios-day0-reset
usr/libexec/mios/mios-db
usr/libexec/mios/mios-delegation-prefilter
usr/libexec/mios/mios-directory-lookup
usr/libexec/mios/mios-discord-send
usr/libexec/mios/mios-discord-status
usr/libexec/mios/mios-docgen
usr/libexec/mios/mios-docs-index
usr/libexec/mios/mios-doctor
usr/libexec/mios/mios-enroll-secure-boot
usr/libexec/mios/mios-env-probe
usr/libexec/mios/mios-everything
usr/libexec/mios/mios-find
usr/libexec/mios/mios-finetune
usr/libexec/mios/mios-finetune-dataset
usr/libexec/mios/mios-finetune-serve
usr/libexec/mios/mios-firecrawl
usr/libexec/mios/mios-flatpak
usr/libexec/mios/mios-flatpak-beta-migrate
usr/libexec/mios/mios-flatpak-icon-sanitize
usr/libexec/mios/mios-flatpak-init
usr/libexec/mios/mios-flatpak-install
usr/libexec/mios/mios-flatpak-overrides-apply
usr/libexec/mios/mios-flatpak-preflight
usr/libexec/mios/mios-forgejo-runner-firstboot.sh
usr/libexec/mios/mios-freeipa-enroll.sh
usr/libexec/mios/mios-gen-role-system
usr/libexec/mios/mios-gpu-passthrough
usr/libexec/mios/mios-gui
usr/libexec/mios/mios-gui-launch
usr/libexec/mios/mios-handoff
usr/libexec/mios/mios-hardcode-lint
usr/libexec/mios/mios-hermes-browser
usr/libexec/mios/mios-hermes-dashboard-auth-stub
usr/libexec/mios/mios-hermes-discord-reactions-patch
usr/libexec/mios/mios-hermes-firstboot
usr/libexec/mios/mios-hermes-init-hook
usr/libexec/mios/mios-hermes-soul-sync
usr/libexec/mios/mios-hermes-tail
usr/libexec/mios/mios-host-launch
usr/libexec/mios/mios-html
usr/libexec/mios/mios-ingest
usr/libexec/mios/mios-installer
usr/libexec/mios/mios-keyring-autounlock
usr/libexec/mios/mios-kg
usr/libexec/mios/mios-knowledge-add
usr/libexec/mios/mios-knowledge-search
usr/libexec/mios/mios-lan-status
usr/libexec/mios/mios-launch
usr/libexec/mios/mios-launcher-daemon
usr/libexec/mios/mios-living-wallpaper
usr/libexec/mios/mios-locate
usr/libexec/mios/mios-login-account
usr/libexec/mios/mios-map
usr/libexec/mios/mios-mcp-enable-tier0.sh
usr/libexec/mios/mios-mcp-server
usr/libexec/mios/mios-md
usr/libexec/mios/mios-micro-llm
usr/libexec/mios/mios-model-router
usr/libexec/mios/mios-new
usr/libexec/mios/mios-open-url
usr/libexec/mios/mios-os-control
usr/libexec/mios/mios-os-recipe
usr/libexec/mios/mios-oscap-gate
usr/libexec/mios/mios-oscontrol-health
usr/libexec/mios/mios-owui-apply-knowledge
usr/libexec/mios/mios-owui-apply-suggestions
usr/libexec/mios/mios-owui-apply-system-prompt
usr/libexec/mios/mios-owui-apply-websearch
usr/libexec/mios/mios-owui-bootstrap-admin
usr/libexec/mios/mios-owui-install-computer-use
usr/libexec/mios/mios-owui-install-pipe
usr/libexec/mios/mios-owui-install-tools
usr/libexec/mios/mios-passport
usr/libexec/mios/mios-pc-control
usr/libexec/mios/mios-pc-vision
usr/libexec/mios/mios-pg-query
usr/libexec/mios/mios-policy-arbiter
usr/libexec/mios/mios-powershell
usr/libexec/mios/mios-ps
usr/libexec/mios/mios-rag
usr/libexec/mios/mios-registry
usr/libexec/mios/mios-remember
usr/libexec/mios/mios-remote
usr/libexec/mios/mios-restart
usr/libexec/mios/mios-sandbox-exec
usr/libexec/mios/mios-scheduled-research
usr/libexec/mios/mios-screenshot
usr/libexec/mios/mios-show-image
usr/libexec/mios/mios-skill-clone
usr/libexec/mios/mios-skills
usr/libexec/mios/mios-smart-resize
usr/libexec/mios/mios-sriov-init
usr/libexec/mios/mios-ssh-dev-cmd
usr/libexec/mios/mios-steamcmd
usr/libexec/mios/mios-stresstest
usr/libexec/mios/mios-suggestion-refresh
usr/libexec/mios/mios-summarize
usr/libexec/mios/mios-swarm-pack-firstboot
usr/libexec/mios/mios-sync-theme
usr/libexec/mios/mios-sync-to-root
usr/libexec/mios/mios-sys-env
usr/libexec/mios/mios-system-status
usr/libexec/mios/mios-sysview
usr/libexec/mios/mios-text-edit
usr/libexec/mios/mios-theme-render
usr/libexec/mios/mios-toml-get
usr/libexec/mios/mios-tool-clone
usr/libexec/mios/mios-tool-search
usr/libexec/mios/mios-ttyd-expose
usr/libexec/mios/mios-ttyd-launch
usr/libexec/mios/mios-verify-launch
usr/libexec/mios/mios-viking
usr/libexec/mios/mios-web-extract
usr/libexec/mios/mios-web-search
usr/libexec/mios/mios-webtools-firstboot.sh
usr/libexec/mios/mios-win-scan
usr/libexec/mios/mios-window
usr/libexec/mios/mios-window-active
usr/libexec/mios/mios-windows
usr/libexec/mios/mios-winget
usr/libexec/mios/mios-wsl-flatpak-export-sync.sh
usr/libexec/mios/mios-wsl-flatpak-heal
usr/libexec/mios/mios-wslg-env-import
usr/share/containers/systemd/mios-adguard.container
usr/share/containers/systemd/mios-ceph.container
usr/share/containers/systemd/mios-cockpit-link.container
usr/share/containers/systemd/mios-cpu-node.container
usr/share/containers/systemd/mios-forge.container
usr/share/containers/systemd/mios-forgejo-runner.container
usr/share/containers/systemd/mios-guacamole.container
usr/share/containers/systemd/mios-guacd.container
usr/share/containers/systemd/mios-k3s.container
usr/share/containers/systemd/mios-llm-heavy-alt.container
usr/share/containers/systemd/mios-llm-heavy.container
usr/share/containers/systemd/mios-llm-light.container
usr/share/containers/systemd/mios-llm-worker@.container
usr/share/containers/systemd/mios-open-webui.container
usr/share/containers/systemd/mios-otelcol.container
usr/share/containers/systemd/mios-pgvector.container
usr/share/containers/systemd/mios-pxe-hub.container
usr/share/containers/systemd/mios-searxng.container
usr/share/containers/systemd/mios-webtools-crawl4ai.container
usr/share/containers/systemd/mios-webtools-firecrawl-api.container
usr/share/containers/systemd/mios-webtools-firecrawl-worker.container
usr/share/containers/systemd/mios-webtools-redis.container
usr/share/doc/mios/adr/0001-two-gate-bake-activation.md
usr/share/doc/mios/adr/0002-mios-sys-shared-base-consolidation.md
usr/share/doc/mios/adr/0003-sbom-not-hardcode.md
usr/share/doc/mios/adr/0004-github-forgejo-equal-publisher.md
usr/share/doc/mios/adr/0005-sovereign-run-off-m-drive.md
usr/share/doc/mios/adr/0006-openai-api-only-ai-contract.md
usr/share/doc/mios/adr/0007-governance-model-laws-adrs-spec.md
usr/share/doc/mios/adr/0008-mios-cat-unified-entry-and-minification.md
usr/share/doc/mios/adr/0009-unified-config-surface.md
usr/share/doc/mios/adr/0010-ssot-as-system-dotfiles.md
usr/share/doc/mios/adr/0011-unified-languages-and-file-patterns.md
audit_all.ps1
run-pipeline.ps1
tools/native/Cargo.toml
usr/share/mios/windows/mios-uia-dump.ps1
Uninstall-MiOS.ps1
powershell/profile.ps1
tools/native/mios-version-check/Cargo.toml
usr/share/doc/mios/archive/absorbed-plans-2026-06/AIOS-GAP-IMPLEMENTATION-PLAN-2026-06-14.md
usr/share/doc/mios/archive/absorbed-plans-2026-06/AIOS-MIOS-MASTER-PLAN-2026-06-14.md
usr/share/doc/mios/archive/absorbed-plans-2026-06/MIOS-FRESH-CONTEXT-EXECUTION-PLAN.md
usr/share/doc/mios/guides/bootstrap_install.md
usr/share/doc/mios/guides/hummingbird-distroless.md
usr/share/doc/mios/guides/inference-consolidation.md
usr/share/doc/mios/guides/memory-tiering.md
usr/share/doc/mios/reference/ANTIFAB-FABRICATION-GAPS.md
usr/share/doc/mios/reference/PIPE-LEAK-RELIABILITY-GAPS.md
usr/share/mios/hermes/config-worker.yaml
usr/share/mios/windows/Set-MiOSWallpaper.ps1
