# usr/lib/systemd/system/systemd-logind.service.d/10-grub2-logind-service.conf
# AI-hint: Configures systemd-logind to ensure the system enters the GRUB2 bootloader menu upon a reboot command, allowing users to select kernels or boot options via the systemd service.
[Service]
Environment=SYSTEMD_REBOOT_TO_BOOT_LOADER_MENU=true
