# usr/lib/systemd/system/NetworkManager.service.d/10-mios-container-gate.conf
# AI-hint: Configures NetworkManager to skip execution within OCI containers via ConditionVirtualization=!container to prevent port conflicts and service overlap in containerized environments.
[Unit]
# 'MiOS': Redundant in OCI containers (handled by Podman/Docker)
ConditionVirtualization=!container
