# usr/lib/systemd/system/ollama.service.d/10-mios-wsl2.conf
[Unit]
# Drop-in retained as a stable extension hook, intentionally empty.
#
# Earlier MiOS versions ran ollama as a host RPM/tarball install at
# /usr/local/bin/ollama, and this drop-in skipped it on WSL2 because
# the host binary wasn't shipped on that shape. Ollama now ships as
# a Quadlet container (/usr/share/containers/systemd/ollama.container,
# image docker.io/ollama/ollama:latest) which runs unmodified on every
# shape -- WSL included; CPU fallback when no NVIDIA CDI is present.
# The previous ConditionPathExists / ConditionVirtualization gates
# would now unconditionally skip the unit on WSL and break first-boot
# model pull, so they are removed.
#
# To force-disable ollama, set [quadlets.enable] ollama = false in
# mios.toml and rebuild, or run `systemctl mask ollama.service` on a
# deployed host.
