# usr/lib/systemd/system/ollama.service.d/10-mios-wsl2.conf
[Unit]
# Skip ollama on WSL2: binary may not be present and CUDA is unavailable.
ConditionVirtualization=!wsl
ConditionPathExists=/usr/local/bin/ollama
