# usr/lib/systemd/system/ollama.service.d/20-mios-keepalive.conf
[Service]
# /usr/lib/systemd/system/ollama.service.d/20-mios-keepalive.conf
#
# OLLAMA_KEEP_ALIVE=-1 -- every loaded model stays resident in
# CPU memory / VRAM FOREVER unless an explicit `ollama stop <model>`
# unloads it. Operator directive 2026-05-17: "keep_alive should be
# TRUE for a TRULY Agentic OS--MiOS!" + "the MiOS-Agent themselves
# -- the prompt enhancer / MiOS AI Assistant / dispatcher for all
# MiOS services ... should be keep alive".
#
# Default Ollama behaviour is to unload models after 5 minutes of
# idle (OLLAMA_KEEP_ALIVE=5m). For an always-on agentic stack we
# want the OPPOSITE: every model loaded once stays loaded so the
# next chat / refine / classify / gate call doesn't pay the cold-
# load cost (which on this WSL2 testbed is 3+ minutes per 4 GB
# model). Combined with firstboot's pre-warm, the boot sequence
# loads every tier model and they remain warm for the host's
# uptime.
#
# Per-request keep_alive in the API body still overrides this
# default (mios-delegation-prefilter + mios-daemon both pass
# keep_alive=-1 explicitly), but setting the SERVICE-LEVEL
# default catches every other caller (hermes-agent's chat
# completions, OWUI's auto-titler, ad-hoc curl probes, etc.)
# so the always-on invariant is global rather than per-script.
Environment=OLLAMA_KEEP_ALIVE=-1
