# usr/share/mios/windows/quadlets/mios-ollama-cpu.container
# AI-hint: Quadlet configuration for the mios-ollama-cpu container, providing a CPU-only inference lane on port 11435 for micro-LLMs and embeddings while sharing the model store with the GPU instance.
# AI-related: mios-ollama-cpu, mios-ollama, network-online.target, multi-user.target, default.target, localhost:11435
# mios-ollama-cpu.container -- Windows AI Node variant
# CPU light-lane: micro-LLMs + embeddings.  Port :11435.
# Mirrors usr/share/containers/systemd/mios-ollama-cpu.container.
# Shares the model store with mios-ollama so pulled models are visible
# to both instances without duplication.

[Unit]
Description=MiOS Ollama CPU Light Lane (Windows AI Node)
After=network-online.target

[Container]
Image=docker.io/ollama/ollama:latest
ContainerName=mios-ollama-cpu
# Host networking -- avoids the WSL2 port-forwarder quirk where a bridge
# + PublishPort binding never reliably surfaces on the host interface.
# Matches the mios-ollama-cpu pattern on MiOS-DEV.
Network=host
AutoUpdate=registry

# CPU-only -- no AddDevice.  Keeps the dGPU fully available for mios-ollama.
Environment=OLLAMA_HOST=0.0.0.0:11435
Environment=OLLAMA_ORIGINS=*
# Shared model store with the GPU instance.
Environment=OLLAMA_MODELS=/var/lib/ollama/models
Environment=HOME=/var/lib/ollama-cpu
# 2026-06-01j runaway fix: 4 concurrent CPU gens each peg all cores -> load 387.
# 2 = one heavy agent gen + one slot for the fast router/refine micro-LLMs.
Environment=OLLAMA_NUM_PARALLEL=2
Environment=OLLAMA_MAX_LOADED_MODELS=2

Volume=/var/lib/ollama:/var/lib/ollama:Z
Volume=/var/lib/ollama-cpu:/var/lib/ollama-cpu:Z

Label=org.opencontainers.image.title=mios-ollama-cpu
Label=org.opencontainers.image.url=http://localhost:11435/

[Service]
Restart=on-failure
RestartSec=10
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
