# usr/lib/systemd/system/systemd-journald.service.d/nice.conf
# AI-hint: Configures the systemd-journald service with a priority of Nice=-1 to prevent the watchdog from triggering during soft lockups, ensuring stable logging during high-load system states.
[Service]
# Set Nice=-1 to doge watchdog on soft lockups (LP: #1696970).
Nice=-1
