# etc/containers/systemd/mios-ai.container
# /etc/containers/systemd/mios-ai.container
# OpenAI-compatible inference backend.
#
# Defaults policy: enabled by default. ConditionPathExists guards the
# config dir; if /etc/mios/ai/ doesn't exist (incomplete bootstrap),
# the unit no-ops at pre-boot rather than crash-looping.

[Unit]
Description='MiOS' AI inference backend (OpenAI-compatible)
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/etc/mios/ai

[Container]
# Image / Port / User-Group placeholders -- resolved at image build
# time by automation/15-render-quadlets.sh from mios.toml
# [image.sidecars].localai, [ports].localai, [services.localai].uid,
# [services.localai].gid (via tools/lib/userenv.sh's MIOS_* env vars).
# Editing mios.toml is the only place these values change; the
# rendered file in the deployed image carries the resolved literals.
Image=docker.io/${MIOS_LOCALAI_IMAGE:-localai/localai:v2.20.0}
ContainerName=mios-ai
# MiOS internal networks: mios.network is the historical core bridge
# (mios-forge / mios-cockpit-link / mios-mcp etc.); ai-net.network is
# the dedicated AI-stack bridge that mios-hermes / mios-webui /
# mios-ollama / mios-searxng all share so the gateway can address
# them by container name without leaving the bridge. Multi-network
# membership keeps mios-ai reachable from BOTH worlds.
Network=${MIOS_QUADLET_NETWORK:-mios.network}
Network=ai-net.network
PublishPort=${MIOS_PORT_LOCALAI:-8080}:8080
# LocalAI v4 layout: /build/models is canonical (matches MODELS_PATH /
# upstream Containerfile WORKDIR), but /data/{outputs,collections} is
# hard-coded for image+audio generation and the AgentJobService vector
# store. Earlier MiOS Quadlets remapped these via LOCALAI_*_PATH env
# vars, but v4.1.3 ignores those overrides and falls back to /data/...
# anyway -- yielding `mkdir /data/outputs: permission denied` because
# /data is created at root:root inside the image and mios-ai runs as
# UID 817. Mount the host-side dirs straight onto the canonical paths
# so the binary writes where it expects, with the correct ownership.
Volume=/srv/ai/models:/build/models:Z
Volume=/srv/ai/mcp:/build/mcp:Z
Volume=/srv/ai/outputs:/data/outputs:Z
Volume=/srv/ai/collections:/data/collections:Z
Volume=/etc/mios/ai:/etc/mios/ai:Z,ro
Environment=MODELS_PATH=/build/models
Environment=THREADS=4
Environment=CONTEXT_SIZE=4096
# Bind 0.0.0.0 so the OpenAI-API endpoint is reachable from the host
# LAN, sibling containers on mios.network, and the loopback when MiOS
# is deployed as a VM, container, or bare-metal host. The PublishPort
# mapping on the host side is already 0.0.0.0:${MIOS_PORT_LOCALAI} by
# Podman default.
Environment=ADDRESS=0.0.0.0:8080
Environment=API_KEY=
Environment=CORS=true
# Use NUMERIC UID/GID -- the LocalAI image has no `mios-ai` user in its
# /etc/passwd, so a name-based User= lookup fails. UID/GID resolved
# from mios.toml [services.localai].uid/.gid (default 817) and pinned
# in usr/lib/sysusers.d/50-mios-ai.conf so the host's mios-ai user has
# the same UID and shared-volume ownership stays consistent.
User=${MIOS_LOCALAI_UID:-817}
Group=${MIOS_LOCALAI_GID:-817}

# OCI image labels for Podman Desktop. Values MUST NOT contain spaces
# -- Quadlet's Label= parser word-splits on whitespace and emits each
# word as a separate --label arg, breaking the entire podman-run
# command line.
Label=org.opencontainers.image.title=mios-ai
Label=org.opencontainers.image.url=http://localhost:8080/v1/models
Label=org.opencontainers.image.documentation=https://platform.openai.com/docs/api-reference
Label=io.podman_desktop.openInBrowser=http://localhost:8080/v1/models

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=900s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
