# etc/containers/systemd/mios-ai.container
# /etc/containers/systemd/mios-ai.container
# OpenAI-compatible inference backend.
#
# Defaults policy: enabled by default. ConditionPathExists guards the
# config dir; if /etc/mios/ai/ doesn't exist (incomplete bootstrap),
# the unit no-ops at pre-boot rather than crash-looping.

[Unit]
Description='MiOS' AI inference backend (OpenAI-compatible)
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/etc/mios/ai

[Container]
Image=docker.io/localai/localai:latest
ContainerName=mios-ai
PublishPort=8080:8080
Volume=/srv/ai/models:/build/models:Z
Volume=/srv/ai/mcp:/build/mcp:Z
Volume=/srv/ai/outputs:/build/outputs:Z
Volume=/etc/mios/ai:/etc/mios/ai:Z,ro
Environment=MODELS_PATH=/build/models
# LocalAI image hard-codes /data/outputs as the audio/image generation
# scratch path. The image creates that dir at root:root (rwx by root
# only); mios-ai runs as UID 817 -> EACCES on first write. Redirect to
# the mounted /build/outputs volume which tmpfiles owns to UID 817.
Environment=LOCALAI_IMAGE_PATH=/build/outputs
Environment=LOCALAI_AUDIO_PATH=/build/outputs
Environment=LOCALAI_UPLOAD_PATH=/build/outputs
Environment=THREADS=4
Environment=CONTEXT_SIZE=4096
# Bind 0.0.0.0 so the OpenAI-API endpoint is reachable from the host
# LAN, sibling containers on mios.network, and the loopback when MiOS
# is deployed as a VM, container, or bare-metal host. The PublishPort
# mapping on the host side is already 0.0.0.0:8080 by Podman default.
Environment=ADDRESS=0.0.0.0:8080
Environment=API_KEY=
Environment=CORS=true
# Use NUMERIC UID/GID -- the LocalAI image has no `mios-ai` user in its
# /etc/passwd, so a name-based User= lookup fails with "unable to find
# user mios-ai: no matching entries in passwd file". UID 817 is pinned
# in usr/lib/sysusers.d/50-mios-ai.conf so the host's mios-ai user has
# the same UID, and shared-volume ownership stays consistent.
User=817
Group=817

# OCI image labels for Podman Desktop. Values MUST NOT contain spaces
# -- Quadlet's Label= parser word-splits on whitespace and emits each
# word as a separate --label arg, breaking the entire podman-run
# command line.
Label=org.opencontainers.image.title=mios-ai
Label=org.opencontainers.image.url=http://localhost:8080/v1/models
Label=org.opencontainers.image.documentation=https://platform.openai.com/docs/api-reference
Label=io.podman_desktop.openInBrowser=http://localhost:8080/v1/models

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=900s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
