# etc/containers/systemd/mios-ai.container
# /etc/containers/systemd/mios-ai.container
# OpenAI-compatible inference backend.
#
# Defaults policy: enabled by default. ConditionPathExists guards the
# config dir; if /etc/mios/ai/ doesn't exist (incomplete bootstrap),
# the unit no-ops at pre-boot rather than crash-looping.

[Unit]
Description='MiOS' AI inference backend (OpenAI-compatible)
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/etc/mios/ai

[Container]
Image=docker.io/localai/localai:latest
ContainerName=mios-ai
PublishPort=8080:8080
Volume=/srv/ai/models:/build/models:Z
Volume=/srv/ai/mcp:/build/mcp:Z
Volume=/etc/mios/ai:/etc/mios/ai:Z,ro
Environment=MODELS_PATH=/build/models
Environment=THREADS=4
Environment=CONTEXT_SIZE=4096
# Bind 0.0.0.0 so the OpenAI-API endpoint is reachable from the host
# LAN, sibling containers on mios.network, and the loopback when MiOS
# is deployed as a VM, container, or bare-metal host. The PublishPort
# mapping on the host side is already 0.0.0.0:8080 by Podman default.
Environment=ADDRESS=0.0.0.0:8080
Environment=API_KEY=
Environment=CORS=true
User=mios-ai
Group=mios-ai

# OCI image labels for Podman Desktop. Architectural Law 5: this is the
# single MIOS_AI_ENDPOINT every agent on the system targets; surfacing
# it in Podman Desktop lets operators sanity-check that the OpenAI-API-
# shaped surface is reachable without leaving the GUI.
Label=org.opencontainers.image.title=MiOS AI inference (LocalAI)
Label=org.opencontainers.image.description=OpenAI-API-compatible /v1 endpoint at port 8080 (chat/completions, responses, embeddings, models).
Label=org.opencontainers.image.url=http://localhost:8080/v1/models
Label=org.opencontainers.image.documentation=https://platform.openai.com/docs/api-reference
Label=io.podman_desktop.openInBrowser=http://localhost:8080/v1/models

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=900s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
