# etc/containers/systemd/mios-ceph.container
# /etc/containers/systemd/mios-ceph.container
# Minimal Ceph monitor for storage orchestration.
#
# Image tag: bootstrap overrides this with the site-specific Ceph version.
# v18 (Reef) is the 'MiOS' default stable branch tag -- less mutable than :latest.
# Ceph and K3s run as root (uid 0) by necessity: Ceph requires privileged
# storage access; both are documented architectural exceptions to UNPRIVILEGED-EXECUTION.

[Unit]
Description='MiOS' Ceph Monitor (Podman-native)
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips on hosts that aren't
# Ceph-configured. Ceph requires a configured cluster; without
# /etc/ceph/ceph.conf the monitor cannot bootstrap, so skip cleanly at
# pre-boot. Nested containers cannot run a Ceph mon.
ConditionPathExists=/etc/ceph/ceph.conf
ConditionVirtualization=!container

[Container]
Image=quay.io/ceph/ceph:latest
ContainerName=mios-ceph
Network=mios.network
Volume=/var/lib/ceph:/var/lib/ceph:Z
Volume=/etc/ceph:/etc/ceph:Z
Exec=mon
# Documented exception (see header): User=/Group= declared explicitly as
# root per LAW 6, even though Ceph requires uid 0.
User=root
Group=root

[Service]
Restart=always
Delegate=yes

[Install]
WantedBy=multi-user.target
