# etc/containers/systemd/mios-forge.container
# /etc/containers/systemd/mios-forge.container
# 'MiOS' self-hosted Git forge (Forgejo upstream).
#
# Rationale: Forgejo is the resource-fit-for-MiOS forge per the analysis
# in usr/share/doc/mios/reference/credits.md / usr/share/doc/mios/reference/sources.md (idle 150-300 MB, single-core sufficient for
# small teams, GPLv3+ governance under Codeberg e.V., GitHub Actions-
# compatible CI engine via the Forgejo Runner so existing
# .github/workflows/ runs unmodified). The image is pulled from
# codeberg.org (Forgejo's upstream registry) by default.
#
# Defaults policy: enabled by default. ConditionPathIsDirectory guards
# the config dir; if /etc/mios/forge/ doesn't exist (incomplete bootstrap)
# the unit no-ops at pre-boot rather than crash-looping. !container also
# skips inside nested containers where the bind-mount pattern would clash.

[Unit]
Description='MiOS' Git forge (Forgejo)
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/etc/mios/forge
# Forgejo IS the self-replication loop's git origin -- must run on every
# shape including WSL. systemd reports WSL as a `container` virt variant,
# so a bare `!container` would skip WSL and break the entire push -> CI
# -> bootc switch loop on Windows hosts.
ConditionVirtualization=|!container
ConditionVirtualization=|wsl

[Container]
# Image / Port / UID placeholders -- resolved at image build time by
# automation/15-render-quadlets.sh from mios.toml [image.sidecars].forge,
# [ports].forge_http, [ports].forge_ssh, [services.forge].uid/.gid.
Image=${MIOS_FORGE_IMAGE:-codeberg.org/forgejo/forgejo:12}
ContainerName=mios-forge
Network=${MIOS_QUADLET_NETWORK:-mios.network}
# Also join ai-net so mios-code-server (and any future editor / agent
# container) can reach the forge by container DNS as mios-forge:3000.
Network=ai-net.network
# forge_http = web UI (reverse-proxy via Cockpit/Caddy if desired).
# forge_ssh  = git+ssh, mapped to container :22; non-22 host port so sshd
# (also on :22) and the forge's git transport coexist on the same host.
PublishPort=${MIOS_PORT_FORGE_HTTP:-3000}:3000
PublishPort=${MIOS_PORT_FORGE_SSH:-2222}:22

# OCI image labels. Values MUST NOT contain whitespace -- Quadlet's
# Label= parser word-splits on spaces and emits each word as a separate
# --label arg to `podman run`, which then mis-positions the image ref
# and the container fails to start. Use hyphens/dots/no-break-space
# encoding if a multi-word description is genuinely needed.
Label=org.opencontainers.image.title=mios-forge
Label=org.opencontainers.image.url=http://localhost:3000/
Label=org.opencontainers.image.documentation=https://forgejo.org/docs/latest/
Label=io.podman_desktop.openInBrowser=http://localhost:3000/

# Persistent state: the Forgejo image expects /data to hold app state +
# the SQLite DB + git repository bytes. /srv/mios/forge is the FHS-3.0
# 'data served by the system' location declared in usr/lib/tmpfiles.d/
# mios-forge.conf and chowned mios-forge:mios-forge there.
Volume=/srv/mios/forge:/data:Z
Volume=/etc/mios/forge:/etc/mios/forge:Z,ro

# Identity inside the container. The upstream image creates a 'git' user
# with USER_UID/USER_GID env vars at startup; we pin those to MiOS's
# 50-mios-services.conf reservation (816 = next free slot in 810-829).
Environment=USER=git
Environment=USER_UID=${MIOS_FORGE_UID:-816}
Environment=USER_GID=${MIOS_FORGE_GID:-816}

# Forgejo configuration (env-var-driven; no app.ini needed for default
# deployment). Format: FORGEJO__<section>__<key>=<value>. Operators
# wanting richer config can drop a custom app.ini under /etc/mios/forge/
# and add Environment=FORGEJO_APP_INI_PATH=/etc/mios/forge/app.ini.
Environment=FORGEJO__server__ROOT_URL=http://localhost:${MIOS_PORT_FORGE_HTTP:-3000}/
# 127.0.0.1, NOT 0.0.0.0. gitea/forgejo's net.Listen("tcp",
# "0.0.0.0:3000") path internally upgrades to AF_INET6 dual-stack
# (`*:3000` in ss output, AF_INET6-with-V6ONLY=0). WSL2 NAT-mode
# localhostForwarding only catches AF_INET binds; AF_INET6 dual-
# stack listens silently miss the forwarder and Windows times out
# on 127.0.0.1:3000. Forcing HTTP_ADDR=127.0.0.1 makes gitea bind
# AF_INET on 127.0.0.1:3000 explicitly, which WSL2 forwards.
# Operator-flagged 2026-05-11 (`localhost:3000 timed out` after
# searxng+cockpit IPv4 fixes were in).
# Trade-off: LAN-side reach (phone, second laptop) needs an extra
# socat hop or Windows-side netsh portproxy from the host's LAN IP
# to the dev VM's NAT IP. Tracked separately; Windows-localhost
# access is the higher-volume use case so that wins for the default.
Environment=FORGEJO__server__HTTP_ADDR=127.0.0.1
Environment=GITEA__server__HTTP_ADDR=127.0.0.1
Environment=FORGEJO__server__HTTP_PORT=${MIOS_PORT_FORGE_HTTP:-3000}
Environment=FORGEJO__server__SSH_PORT=${MIOS_PORT_FORGE_SSH:-2222}
Environment=FORGEJO__server__SSH_LISTEN_HOST=0.0.0.0
Environment=FORGEJO__server__SSH_LISTEN_PORT=22
Environment=FORGEJO__database__DB_TYPE=sqlite3
Environment=FORGEJO__database__PATH=/data/forgejo.db
Environment=FORGEJO__service__DISABLE_REGISTRATION=true
Environment=FORGEJO__service__REQUIRE_SIGNIN_VIEW=false
Environment=FORGEJO__security__INSTALL_LOCK=true
Environment=FORGEJO__log__LEVEL=Info
Environment=FORGEJO__log__MODE=console
# Architectural Law 5: nothing in the forge calls vendor AI URLs.
# Forgejo has no built-in AI surface; this comment exists so postcheck
# #12 has a documented reason for the absence rather than an oversight.

# DOCUMENTED ARCHITECTURAL EXCEPTION (UNPRIVILEGED-QUADLETS, LAW 6):
# The upstream Forgejo image is built around s6-overlay; the entrypoint
# MUST start as root so it can (a) `sed` /etc/group + /etc/passwd to
# pin the `git` user to USER_UID/USER_GID (816), (b) chown /data and
# the runtime dirs to that uid, and (c) start s6-svscan, which writes
# its supervisor lock under /run. All three steps fail with
# "Permission denied" the moment the container starts as a non-root
# uid -- the symptoms observed during the 2026-05-05 bootstrap regression:
#     sed: can't create temp file '/etc/groupXXXXXX': Permission denied
#     sed: can't create temp file '/etc/passwdXXXXXX': Permission denied
#     s6-svscan: fatal: unable to open .s6-svscan/lock: Permission denied
# After the entrypoint finishes its root-only setup, s6-overlay drops
# to the `git` user (uid 816) for the actual long-lived forgejo daemon
# process, so the running surface is still unprivileged. Host-side
# ownership of /srv/mios/forge is set to uid 816 by usr/lib/tmpfiles.d/
# mios-forge.conf. LAW 6 (UNPRIVILEGED-QUADLETS) is satisfied at the
# runtime layer that matters (the long-lived forgejo process), not
# the millisecond root-bootstrap window the image's design requires.
# Joins the existing exception list: mios-ceph, mios-k3s, and
# mios-forgejo-runner -- all four declare User= explicitly so the
# 99-postcheck.sh "every Quadlet declares User=" gate passes; the
# docs/architecture.md row notes "starts root, drops privileges via
# upstream image entrypoint" as the load-bearing distinction.
User=0
Group=0

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
