# etc/containers/systemd/mios-k3s.container
# /etc/containers/systemd/mios-k3s.container
# Rootless K3s orchestration via Podman.
#
# Privileged=true: K3s requires full host networking, kernel namespaces, and
# eBPF -- specific capability grants cannot replace these. This is a documented
# architectural exception to UNPRIVILEGED-EXECUTION (see ARCHITECTURE.md).
#
# User= omitted: K3s runs as root by necessity (writes to /run/k3s, manages
# cgroups, kernel netfilter). This is a second documented exception.
#
# K3S_TOKEN: This is the 'MiOS' default cluster token. The bootstrap repository
# overrides it via a systemd drop-in at ignition time:
#   /etc/containers/systemd/mios-k3s.container.d/token.conf
# NEVER change this default to a real secret -- it is visible in the image.

[Unit]
Description='MiOS' K3s Service (Podman-native)
After=network-online.target
Wants=network-online.target
# Defaults policy: enabled by default; auto-skips on incompatible hosts.
# K3s needs cgroup v2, eBPF, and full kernel namespaces -- none of which
# work reliably under WSL2 or in nested containers.
ConditionVirtualization=!wsl
ConditionVirtualization=!container

[Container]
Image=docker.io/rancher/k3s:latest
ContainerName=mios-k3s
Privileged=true
Environment=K3S_TOKEN=mios-cluster-secret
Volume=/var/lib/rancher/k3s:/var/lib/rancher/k3s:Z
Volume=/usr/share/mios/k3s-manifests:/var/lib/rancher/k3s/server/manifests:ro
Exec=server --disable=traefik
# Documented exception (see header): User=/Group= declared explicitly as
# root per LAW 6, even though K3s requires Privileged=true and uid 0.
User=root
Group=root

[Service]
Restart=always
Delegate=yes

[Install]
WantedBy=multi-user.target
