# etc/containers/systemd/mios-ai.container
[Unit]
Description=MiOS AI Gateway (LocalAI)
After=network-online.target

[Container]
Image=quay.io/go-skynet/local-ai:latest
Environment=MODELS_PATH=/v1/models
Volume=/v1:/v1:ro
Volume=/usr/share/mios/ai:/usr/share/mios/ai:ro
PublishPort=8080:8080
Label=io.mios.role=ai-gateway
Network=mios.network

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
