# etc/containers/systemd/mios-code-server.container
# /etc/containers/systemd/mios-code-server.container
# code-server -- VS Code in a browser. Pairs with mios-forge for a
# full editor + git workflow over the local LAN without leaving the
# browser.
#
# Upstream: https://github.com/coder/code-server
#
# Operator workflow:
#   1. Open http://localhost:8800/ in any browser.
#   2. Log in with the global MiOS password (mios.toml [identity].
#      default_password; default "mios").
#   3. Inside code-server's terminal:
#        git clone http://mios-forge:3000/<user>/<repo>.git
#      The container reaches mios-forge via ai-net.network DNS, so the
#      bridge alias works without leaking host loopback into the
#      sandbox. Forgejo HTTP basic auth uses the same mios/mios pair
#      (forge admin password is also MIOS_DEFAULT_PASSWORD).

[Unit]
Description='MiOS' code-server (VS Code in a browser)
Documentation=https://github.com/coder/code-server
After=network-online.target mios-forge.service
Wants=network-online.target

[Container]
# Image / port placeholders rendered at image build time by
# automation/15-render-quadlets.sh from mios.toml [image.sidecars].
# code_server + [ports].code_server.
Image=${MIOS_CODE_SERVER_IMAGE:-ghcr.io/coder/code-server:latest}
ContainerName=mios-code-server
# Run as the mios operator. MIOS_CODE_SERVER_UID / MIOS_CODE_SERVER_GID
# are resolved at image-build time by 15-render-quadlets.sh from
# `id -u mios` / `id -g mios` on the host so the bind-mount ownership
# below lines up regardless of which uid the base distro auto-allocated
# for the mios sysuser (podman-machine-os 6.0 happens to land mios at
# 992, not the 1000 the MiOS overlay assumes elsewhere).
# Run as the in-container 'coder' user (uid 1000 in the upstream image).
# The Volume= mount below points at a dedicated host dir chown'd
# 1000:1000 by tmpfiles.d so the perms line up regardless of which uid
# the dev VM's base distro auto-allocated for mios (podman-machine-os
# 6.0 lands mios at 992, breaking the host-home-mount pattern).
User=1000
Group=1000
# Network=host: bind on the dev VM netns directly so localhost:<port>
# hits the editor without netavark's port-forwarder (which intermittently
# times out on WSL2/podman-machine-os 6.0). Forge is reachable as
# http://localhost:3000/ from code-server's terminal in this mode.
Network=host
SecurityLabelDisable=true

# Dedicated workspace dir (NOT /var/home/mios) so the in-container uid
# 1000 matches the host dir owner. Operator can symlink /var/home/mios/
# git into /var/lib/mios/code-server to surface their repos in the
# editor sidebar without UID-drift gymnastics.
Volume=/var/lib/mios/code-server:/home/coder

# Auth: MIOS_DEFAULT_PASSWORD baked at image-build time from
# mios.toml [identity].default_password via 15-render-quadlets.sh.
# No runtime EnvironmentFile= -- Quadlet [Container] doesn't honor
# the systemd `-` optional-file prefix and concatenates `-/etc/...`
# into the path arg, breaking podman --env-file. Rotate by editing
# mios.toml + rebuilding (or hand-edit this Environment= line on a
# live host then `systemctl daemon-reload && systemctl restart`).
Environment=PASSWORD=${MIOS_DEFAULT_PASSWORD:-mios}
Environment=SUDO_PASSWORD=${MIOS_DEFAULT_PASSWORD:-mios}
Environment=DEFAULT_WORKSPACE=/home/coder
# The upstream image's ENTRYPOINT hardcodes --bind-addr 0.0.0.0:8080.
# We use port 8080 with Network=host so the container's listener IS
# the host's. BIND_ADDR is kept for completeness but the entrypoint
# arg wins.
Environment=BIND_ADDR=0.0.0.0:${MIOS_PORT_CODE_SERVER:-8080}
AutoUpdate=registry

Label=org.opencontainers.image.title=mios-code-server
Label=org.opencontainers.image.url=http://localhost:${MIOS_PORT_CODE_SERVER:-8080}/
Label=org.opencontainers.image.documentation=https://github.com/coder/code-server
Label=io.podman_desktop.openInBrowser=http://localhost:${MIOS_PORT_CODE_SERVER:-8080}/

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=600s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
