# etc/containers/systemd/mios-hermes-workspace.container
# /etc/containers/systemd/mios-hermes-workspace.container
# Hermes Workspace -- chat front-end for the Hermes-Agent stack.
#
# Chain: ollama -> mios-hermes -> mios-hermes-workspace (default chat UI
# the operator opens at http://localhost:3030/).
#
# Upstream: https://github.com/outsourc-e/hermes-workspace
#
# Auth: HERMES_API_TOKEN MUST equal API_SERVER_KEY from
# /etc/mios/hermes/api.env (generated by mios-hermes-firstboot.service).
# HERMES_PASSWORD is the workspace session login, generated by
# mios-hermes-workspace-firstboot.service into
# /etc/mios/hermes-workspace/workspace.env.

[Unit]
Description='MiOS' Hermes Workspace (web frontend for Hermes Agent)
After=network-online.target mios-hermes.service
Wants=network-online.target mios-hermes.service
ConditionPathExists=/etc/mios/hermes/api.env

[Container]
# Image / port placeholders rendered at image build time by
# automation/15-render-quadlets.sh from mios.toml.
Image=ghcr.io/outsourc-e/hermes-workspace:latest
ContainerName=mios-hermes-workspace
# Upstream node:22-slim image declares no USER and expects root to
# chown /app/data + spawn /scripts/pty-helper. Law 6 exception
# matches mios-ceph and mios-k3s.
User=0
Group=0
Network=mios.network
Network=ai-net.network
PublishPort=3030:3030
AutoUpdate=registry

# HERMES_API_TOKEN sourced from /etc/mios/hermes/api.env. Quadlet
# Environment= directives don't expand ${VAR} from EnvironmentFile=
# contents, so the env file is loaded directly here.
EnvironmentFile=/etc/mios/hermes/api.env
EnvironmentFile=/etc/mios/hermes-workspace/workspace.env
Environment=HERMES_API_URL=http://mios-hermes:8642
# Override upstream PORT=3000 to keep 3030 (mios.toml [ports].hermes_workspace).
Environment=PORT=3030
# Browsers drop Secure cookies over plain HTTP; the workspace's
# NODE_ENV=production default sets Secure. Disable so first-time login
# works without a TLS proxy in front. Flip back to 1 in
# /etc/mios/hermes-workspace/workspace.env when fronting via HTTPS.
Environment=COOKIE_SECURE=0

Label=org.opencontainers.image.title=mios-hermes-workspace
Label=org.opencontainers.image.url=http://localhost:3030/
Label=org.opencontainers.image.documentation=https://github.com/outsourc-e/hermes-workspace
Label=io.podman_desktop.openInBrowser=http://localhost:3030/

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=600s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
