# etc/containers/systemd/mios-cockpit-link.container
# /etc/containers/systemd/mios-cockpit-link.container
# Podman Desktop UI surfacing for the host-side Cockpit web console.
#
# Cockpit on 'MiOS' runs as a host systemd service (cockpit.socket bound
# to :9090), not in a container -- so Podman Desktop's container view
# normally has no way to render a clickable link to it. This Quadlet is
# a discovery shim: a tiny socat forwarder that publishes
#
#   localhost:19090   ->   host:9090   (Cockpit web console)
#
# When Podman Desktop scans running containers it sees the published
# port and the OCI labels below, and renders a clickable "open" link
# straight into Cockpit. The alternate port (19090) is required because
# the host's :9090 is already bound by cockpit-ws and a second bind
# would collide.
#
# Footprint: alpine/socat is ~3MB and the container is essentially idle
# (single TCP-relay loop). Disable it via /etc/mios/profile.toml's
# [quadlets.enable] table if you don't want this discovery surface.

[Unit]
Description='MiOS' Cockpit web-console discovery shim (Podman Desktop link)
Documentation=https://cockpit-project.org/
After=network-online.target cockpit.socket
Wants=network-online.target
# Skip when cockpit isn't installed or hasn't bound :9090 yet -- no point
# in advertising a dead URL.
ConditionPathExists=/usr/lib/systemd/system/cockpit.socket
# Skip on WSL2. The discovery shim's whole purpose is to give the
# Linux-host Podman Desktop a clickable link to the host's Cockpit
# console; on WSL2 deployments Podman Desktop runs on the *Windows*
# side and reaches Cockpit at https://localhost:9090/ directly
# (Windows port-forwards loopback into the WSL distro), so there is
# no value in spinning up an alpine/socat relay -- and on the WSL2
# kernel the systemd-executor reliably trips
# `Failed to spawn executor: Device or resource busy` on this Quadlet
# when the rootless namespace plumbing collides with WSLInterop
# pidfd quirks. Saving the resource churn outweighs the (zero)
# discoverability benefit on this shape.
ConditionVirtualization=!wsl
ConditionVirtualization=!container

[Container]
Image=docker.io/alpine/socat:latest@sha256:44e7b87dfeddb0f999ce75eedfa140dc97ea5bf7decbd04b2cef6adb5cc331fe
ContainerName=mios-cockpit-link
# Single MiOS internal network -- KISS, every Quadlet on the same
# bridge so sibling lookups by container name work without host-
# loopback gymnastics. The host.containers.internal hop below is
# still needed because cockpit-ws runs as a HOST service, not a
# Quadlet -- it doesn't have a sibling-name on this bridge.
Network=mios.network
# Forward 19090 (container, bound 0.0.0.0) -> host:9090 (cockpit).
# host.containers.internal is Podman's documented loopback-to-host
# alias (mirrors Docker's host.docker.internal); falls back to the
# well-known 10.88.0.1 gateway if the alias isn't resolvable.
# bind=0.0.0.0 makes the listener reachable from the host LAN,
# sibling containers, and the loopback alike (matches the access
# surface contract every other MiOS service holds to).
Exec=tcp-listen:${MIOS_PORT_COCKPIT_LINK:-19090},bind=0.0.0.0,fork,reuseaddr tcp:host.containers.internal:${MIOS_PORT_COCKPIT:-9090}
PublishPort=0.0.0.0:${MIOS_PORT_COCKPIT_LINK:-19090}:${MIOS_PORT_COCKPIT_LINK:-19090}

# OCI labels Podman Desktop renders. Values MUST NOT contain spaces --
# Quadlet's Label= parser word-splits on whitespace and breaks the
# generated podman-run command. The url label is the canonical Cockpit
# endpoint (https; cockpit-ws issues a self-signed cert by default and
# clients accept it on first visit).
Label=org.opencontainers.image.title=mios-cockpit-link
Label=org.opencontainers.image.url=https://localhost:9090/
Label=org.opencontainers.image.documentation=https://cockpit-project.org/documentation.html
Label=io.podman_desktop.openInBrowser=https://localhost:9090/

# Architectural Law 6: unprivileged Quadlet with explicit User/Group.
# Use the `nobody` uid/gid (65534). alpine/socat doesn't need any
# host-side file access; the discovery shim is a pure TCP relay between
# container :19090 and host :9090. Numeric resolution avoids depending
# on a username that doesn't exist in alpine's /etc/passwd.
User=65534
Group=65534

[Service]
Restart=on-failure
RestartSec=15s
TimeoutStartSec=60s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
