# etc/containers/systemd/mios-cockpit-link.container
# /etc/containers/systemd/mios-cockpit-link.container
# Podman Desktop UI surfacing for the host-side Cockpit web console.
#
# Cockpit on 'MiOS' runs as a host systemd service (cockpit.socket bound
# to :9090), not in a container -- so Podman Desktop's container view
# normally has no way to render a clickable link to it. This Quadlet is
# a discovery shim: a tiny socat forwarder that publishes
#
#   localhost:19090   ->   host:9090   (Cockpit web console)
#
# When Podman Desktop scans running containers it sees the published
# port and the OCI labels below, and renders a clickable "open" link
# straight into Cockpit. The alternate port (19090) is required because
# the host's :9090 is already bound by cockpit-ws and a second bind
# would collide.
#
# Footprint: alpine/socat is ~3MB and the container is essentially idle
# (single TCP-relay loop). Disable it via /etc/mios/profile.toml's
# [quadlets.enable] table if you don't want this discovery surface.

[Unit]
Description='MiOS' Cockpit web-console discovery shim (Podman Desktop link)
Documentation=https://cockpit-project.org/
After=network-online.target cockpit.socket
Wants=network-online.target
# Skip when cockpit isn't installed or hasn't bound :9090 yet -- no point
# in advertising a dead URL.
ConditionPathExists=/usr/lib/systemd/system/cockpit.socket
# Pure userspace TCP relay -- safe on every shape including WSL where
# Cockpit is also reachable on :9090 inside the distro.
ConditionVirtualization=|!container
ConditionVirtualization=|wsl

[Container]
Image=docker.io/alpine/socat:latest
ContainerName=mios-cockpit-link
# Forward 19090 (container, bound 0.0.0.0) -> host:9090 (cockpit).
# host.containers.internal is Podman's documented loopback-to-host
# alias (mirrors Docker's host.docker.internal); falls back to the
# well-known 10.88.0.1 gateway if the alias isn't resolvable.
# bind=0.0.0.0 makes the listener reachable from the host LAN,
# sibling containers, and the loopback alike (matches the access
# surface contract every other MiOS service holds to).
Exec=tcp-listen:19090,bind=0.0.0.0,fork,reuseaddr tcp:host.containers.internal:9090
PublishPort=0.0.0.0:19090:19090

# OCI labels Podman Desktop renders. Values MUST NOT contain spaces --
# Quadlet's Label= parser word-splits on whitespace and breaks the
# generated podman-run command. The url label is the canonical Cockpit
# endpoint (https; cockpit-ws issues a self-signed cert by default and
# clients accept it on first visit).
Label=org.opencontainers.image.title=mios-cockpit-link
Label=org.opencontainers.image.url=https://localhost:9090/
Label=org.opencontainers.image.documentation=https://cockpit-project.org/documentation.html
Label=io.podman_desktop.openInBrowser=https://localhost:9090/

# Architectural Law 6: unprivileged Quadlet with explicit User/Group.
# Reuse the mios-ai uid (NUMERIC -- the alpine/socat image has no
# `mios-ai` user, and podman looks up names IN the container's
# /etc/passwd. UID 817 is pinned in usr/lib/sysusers.d/50-mios-ai.conf).
User=817
Group=817

[Service]
Restart=on-failure
RestartSec=15s
TimeoutStartSec=60s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
