# etc/containers/systemd/mios-ai.container
# /etc/containers/systemd/mios-ai.container
# OpenAI-compatible inference backend.
#
# Defaults policy: enabled by default. ConditionPathExists guards the
# config dir; if /etc/mios/ai/ doesn't exist (incomplete bootstrap),
# the unit no-ops at pre-boot rather than crash-looping.

[Unit]
Description='MiOS' AI inference backend (OpenAI-compatible)
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/etc/mios/ai

[Container]
Image=docker.io/localai/localai:latest
ContainerName=mios-ai
PublishPort=8080:8080
Volume=/srv/ai/models:/build/models:Z
Volume=/srv/ai/mcp:/build/mcp:Z
Volume=/etc/mios/ai:/etc/mios/ai:Z,ro
Environment=MODELS_PATH=/build/models
Environment=THREADS=4
Environment=CONTEXT_SIZE=4096
User=mios-ai
Group=mios-ai

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=900s
Delegate=yes

[Install]
WantedBy=multi-user.target default.target
