# etc/containers/systemd/mios-open-webui.container
# /etc/containers/systemd/mios-open-webui.container
# 'MiOS' Open WebUI - rich LLM interface.
#
# Open WebUI provides a ChatGPT-like browser UI fronting the MiOS AI
# stack. AI chain (2026-05-13 architecture):
#   ollama container (:11434)  -- raw inference + embeddings
#     -> Hermes-Agent (:8642/v1) -- DIRECT host install (NOT a
#        container; automation/38-hermes-agent.sh + hermes-agent.service),
#        OpenAI gateway + agent layer
#       -> Open WebUI (this container, :3030) -- browser UI
#
# Open WebUI's OPENAI_API_BASE_URL points at Hermes-Agent so the AGENT
# (with tool-calling / sessions / kanban / skills), not raw ollama, is
# the model Open WebUI serves -- Architectural Law 5 (UNIFIED-AI-REDIRECTS).

[Unit]
Description='MiOS' Open WebUI
# hermes-agent.service is the DIRECT host install (the mios-hermes
# CONTAINER Quadlet is shelved in mios.toml [quadlets.enable]). Open
# WebUI Wants= it but doesn't hard-Require= it -- the UI still starts
# and shows a connection error rather than refusing to boot if the
# agent is mid-restart.
After=network-online.target hermes-agent.service
Wants=network-online.target hermes-agent.service
ConditionPathExists=/etc/mios/hermes/api.env

[Container]
# Image / Port / Network placeholders rendered by automation/15-render-quadlets.sh
Image=${MIOS_OPEN_WEBUI_IMAGE:-ghcr.io/open-webui/open-webui:main}
ContainerName=mios-open-webui
Network=ai-net.network
Network=mios.network
PublishPort=${MIOS_PORT_OPEN_WEBUI:-3030}:8080
AutoUpdate=registry

# Persistent state.
Volume=/var/lib/mios/open-webui:/app/backend/data:Z

# API key + server settings, shared with Hermes-Agent. api.env defines
# OPENAI_API_KEY=<API_SERVER_KEY> as a plain scalar -- that IS the
# bearer Open WebUI sends to the Hermes gateway at OPENAI_API_BASE_URL
# below, and EnvironmentFile= injects it directly into the container
# env. Do NOT add OPENAI_API_KEY_FILE pointing at this same file: the
# _FILE variant makes Open WebUI read the ENTIRE multi-line file as the
# key value (comments + every other VAR= line included), which the
# gateway then rejects. Operator-confirmed 2026-05-14 ("Invalid API
# key" / "You do not have permission" on every Open WebUI call).
EnvironmentFile=/etc/mios/hermes/api.env

# Wiring to the MiOS AI stack -- point at Hermes-Agent (Law 5).
# Hermes-Agent is a DIRECT HOST service now, not a container, so the
# old container-DNS name `mios-hermes` no longer resolves. Podman
# injects `host.containers.internal` into every container as the route
# to the host -- that's how a container reaches a host-side service.
# (The MiOS-DEV host-network drop-in overrides this to localhost:8642.)
Environment=OPENAI_API_BASE_URL=http://host.containers.internal:8642/v1

# Disable Open WebUI's direct-Ollama integration -- all inference must
# route through Hermes-Agent (Law 5). The ollama backend URL is still
# provided for completeness but ENABLE_OLLAMA=False keeps it dormant.
Environment=ENABLE_OLLAMA=False
Environment=OLLAMA_BASE_URL=http://host.containers.internal:11434

# General config
Environment=WEBUI_NAME="MiOS AI"
Environment=WEBUI_SECRET_KEY=mios-stable-secret-change-me
Environment=ENABLE_SIGNUP=False
Environment=DEFAULT_USER_ROLE=admin

# Identity. UID 817 (mios-open-webui)
User=${MIOS_OPEN_WEBUI_UID:-817}
Group=${MIOS_OPEN_WEBUI_GID:-817}

Label=org.opencontainers.image.title=mios-open-webui
Label=org.opencontainers.image.url=http://localhost:3030/
Label=io.podman_desktop.openInBrowser=http://localhost:3030/

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
Delegate=yes

[Install]
WantedBy=multi-user.target
