# etc/containers/systemd/mios-ai.container
# /etc/containers/systemd/mios-ai.container
# OpenAI-compatible inference backend.
# Pinned to a specific tag for reproducibility.

[Unit]
Description=MiOS AI inference backend (OpenAI-compatible)
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/localai/localai:v2.20.0
ContainerName=mios-ai
PublishPort=8080:8080
Volume=/srv/ai/models:/build/models:Z
Volume=/srv/ai/mcp:/build/mcp:Z
Volume=/etc/mios/ai:/etc/mios/ai:Z,ro
Environment=MODELS_PATH=/build/models
Environment=THREADS=4
Environment=CONTEXT_SIZE=4096

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=900s

[Install]
WantedBy=multi-user.target default.target
