# usr/lib/systemd/system/NetworkManager-wait-online.service.d/timeout.conf
[Unit]
# WSL2 uses its own networking bridge — NM wait-online never completes
ConditionVirtualization=!wsl

[Service]
TimeoutStartSec=10
