# usr/lib/systemd/system/NetworkManager-wait-online.service.d/timeout.conf
# AI-hint: Configures a 10s timeout for NetworkManager-wait-online.service and disables it in WSL2 environments to prevent boot hangs caused by non-responsive virtualized network bridges.
# AI-related: NetworkManager-wait-online.service
[Unit]
# WSL2 uses its own networking bridge -- NM wait-online never completes
ConditionVirtualization=!wsl

[Service]
TimeoutStartSec=10
