# usr/lib/systemd/system/user@.service.d/10-mios-wsl-runtime.conf
# /usr/lib/systemd/system/user@.service.d/10-mios-wsl-runtime.conf
#
# Seed XDG_RUNTIME_DIR for the systemd user manager so it can start without
# a logind session. MiOS skips systemd-logind under WSL2 (the upstream
# unit's ConditionVirtualization=!wsl), so pam_systemd never gets to call
# org.freedesktop.login1.CreateSession and the standard "logind exports
# XDG_RUNTIME_DIR via PAM" handshake doesn't happen. The user manager then
# bails out with:
#
#     systemd[N]: Trying to run as user instance, but $XDG_RUNTIME_DIR is
#         not set.
#     user@1000.service: Main process exited, code=exited, status=1/FAILURE
#
# This drop-in declares the env var statically as /run/user/%i (the standard
# per-uid path). mios-wsl-runtime-dir.service is responsible for actually
# creating the directory on tmpfs before user@.service activates -- it's a
# Before= ordering on the user-runtime-dir@<uid>.service that user@.service
# already pulls in. The two together mimic what pam_systemd would have done.

[Service]
Environment=XDG_RUNTIME_DIR=/run/user/%i
