# usr/lib/systemd/system/dbus-broker.service.d/10-mios-wsl2.conf
# Skip dbus-broker on WSL2 (WSL kernels strip the audit subsystem, causing broker to crash)
# Standardized on native systemd WSL2 detection (v252+).
[Unit]
ConditionVirtualization=!wsl
