# usr/lib/systemd/system/dbus-broker.service.d/10-mios-no-audit.conf
# dbus-broker upstream ExecStart is `/usr/bin/dbus-broker-launch --scope
# system --audit`. The `--audit` flag wires libaudit into the broker so it
# can log policy decisions to the kernel audit subsystem. WSL2 kernels
# strip CONFIG_AUDIT, so libaudit's netlink socket call returns EPROTONOSUPPORT
# and the broker exited 1 at startup -- which previously cascaded into
# logind/polkit/cockpit "Failed to connect to system scope bus" failures.
#
# Dropping `--audit` makes broker run on every kernel (with-audit and
# without). Audit-flagged hosts (bare-metal Fedora with auditd) lose the
# *D-Bus policy audit trail* but keep all other auditd functionality
# (auditd.service is still enabled in 90-mios.preset). For MiOS the
# trade is worth it: one identical drop-in covers WSL, Hyper-V, QEMU,
# bare-metal, and Podman-WSL deployments without per-shape masking.
[Service]
ExecStart=
ExecStart=/usr/bin/dbus-broker-launch --scope system
