# usr/lib/systemd/system/mios-gpu-amd.service.d/10-virt-gate.conf
[Unit]
# 'MiOS': Skip in containers/WSL2
ConditionVirtualization=!container
ConditionVirtualization=!wsl
