# usr/lib/systemd/system/sshd.service.d/10-mios-wsl2.conf
# /usr/lib/systemd/system/sshd.service.d/10-mios-wsl2.conf
#
# WSL2 with mirrored networking (the MiOS default in mios-bootstrap's
# .wslconfig setup) exposes the Windows host's :22 inside the distro's
# network namespace. When sshd tries to bind 0.0.0.0:22 it collides
# with whatever Windows-side service is on that port (typically
# Microsoft OpenSSH, which most Windows installs ship enabled).
# Result: sshd repeatedly fails with "Bind to port 22 on 0.0.0.0
# failed: Address already in use" and gives up.
#
# This drop-in skips the upstream sshd.service on WSL. Cross-distro
# SSH access on WSL goes through `wsl -d MiOS` from the host or
# through the WSL2 bridged-networking interface. Operators who DO
# need sshd reachable from the LAN on WSL can override this drop-in
# with their own at /etc/systemd/system/sshd.service.d/ that adds
# `Port 2200` (or any free port) to sshd_config.d/ and removes the
# ConditionVirtualization gate.
#
# Bare-metal, Hyper-V, QEMU, and VM-shape deployments are unaffected:
# sshd binds :22 normally because there's no host-network collision.
[Unit]
ConditionVirtualization=!wsl
