# usr/lib/systemd/system/systemd-networkd-wait-online.service.d/10-mios-wsl2.conf
# AI-hint: Configures systemd-networkd-wait-online to skip network availability checks on WSL2 instances to prevent session timeouts during startup.
# Skip network wait on WSL2 to prevent user session timeouts (WSL v0.2.0+)
[Unit]
ConditionVirtualization=!wsl
