# usr/lib/systemd/system/hermes-agent.service.d/20-mios-paths-env.conf
# AI-hint: Defines environment variables for the hermes-agent service, providing absolute paths for Windows binaries (PowerShell, CMD), the Everything search CLI, and the launcher socket for cross-context communication.
# AI-related: /etc/mios/install.env, mios-paths-env, mios-find, mios-windows, mios-apps, mios-pc-control, mios-launcher, hermes-agent.service
# /usr/lib/systemd/system/hermes-agent.service.d/20-mios-paths-env.conf
#
# Operator directive 2026-05-16: "NO HARDCODED PATHS!!! JUST
# ENVIRONMENT ARGS--Hermes-Agent can natively iterate on tools and
# fix or make new skils /tools".
#
# Environment vars consumed by mios-* helpers (mios-find, mios-windows,
# mios-apps, mios-pc-control, etc.) when invoked by hermes. Helpers
# read $EVERYTHING_CLI, $WINDOWS_POWERSHELL, $WINDOWS_CMD, etc. before
# falling back to mios.toml [paths] before falling back to compiled
# defaults. Operators override per-host via a drop-in here OR via
# /etc/mios/install.env (also EnvironmentFile= mounted by the unit).
#
# Adding/changing a path: edit /etc/mios/install.env -> systemctl
# restart hermes-agent. No code edits, no SKILL.md edits, no Modelfile.

[Service]
# Voidtools Everything CLI -- colon-separated probe list. mios-find
# walks left-to-right; first existing path wins.
#
# Quoted with leading " so systemd's Environment= parser preserves
# the literal spaces in "Program Files" / "Program Files (x86)".
# Without quoting, systemd word-splits on the spaces and rejects
# each split fragment as an invalid environment assignment (the
# whole line is dropped, leaving EVERYTHING_CLI unset and mios-find
# falling all the way back to compiled defaults).
Environment="EVERYTHING_CLI=/mnt/m/Programs/Everything/es.exe:/mnt/c/Program Files/Everything/es.exe:/mnt/c/Program Files (x86)/Everything/es.exe:/mnt/c/Tools/Everything/es.exe:/mnt/c/Users/mios/AppData/Local/Programs/Everything/es.exe"

# Windows PowerShell / cmd full paths (systemd user PATH strips
# /mnt/c/Windows/System32, so broker calls need absolute paths).
Environment=WINDOWS_POWERSHELL=/mnt/c/Windows/System32/WindowsPowerShell/v1.0/powershell.exe
Environment=WINDOWS_CMD=/mnt/c/Windows/System32/cmd.exe

# Operator-side broker socket (cross-user dispatch from agent context
# to operators interactive Windows session).
Environment=LAUNCHER_SOCKET=/run/mios-launcher/launcher.sock
