# usr/lib/systemd/system/hermes-agent.service.d/30-quick-stop.conf
# /usr/lib/systemd/system/hermes-agent.service.d/30-quick-stop.conf
#
# Tighten TimeoutStopSec on hermes-agent.service so `mios-restart hermes`
# (and the implicit restart during `mios-ai-reset`, skill re-syncs, etc.)
# completes in ~30 s instead of waiting up to 4 minutes for an in-flight
# session shutdown.
#
# Default behaviour (operator-confirmed regression 2026-05-15): hermes
# inherits TimeoutStopSec=4min from the systemd defaults and hangs in
# 'deactivating' for the full window when SIGTERM lands mid-tool-call.
# With the systemwide `TimeoutStopFailureMode=abort` drop-in in place,
# systemd then SIGABRTs the process at the 4-minute mark -- net effect:
# every restart blocks the operator's terminal for ~4 minutes when the
# gateway is mid-anything.
#
# 30 s is a deliberate tradeoff:
#   * fast enough that quick test loops aren't punished by routine
#     restarts (skill edits, SOUL.md tweaks, mios-ai-reset cycles)
#   * long enough that small in-flight tool calls (<5 s typical)
#     finish naturally before SIGTERM lands
#   * truncates the ~minutes-long calls (large delegations, mios build
#     wrappers) -- those should be backgrounded anyway per SOUL.md
#     rule 10
#
# Operators with a long-running synchronous task should wait for it
# to complete BEFORE restarting; the 30 s cap is for the common case.

[Service]
TimeoutStopSec=30s
