# usr/lib/systemd/system/hermes-agent.service.d/30-quick-stop.conf
# /usr/lib/systemd/system/hermes-agent.service.d/30-quick-stop.conf
#
# Tighten TimeoutStopSec on hermes-agent.service so `mios-restart hermes`
# (and the implicit restart during `mios-ai-reset`, skill re-syncs, etc.)
# completes in ~30 s instead of waiting up to 4 minutes for an in-flight
# session shutdown.
#
# Default behaviour (operator-confirmed regression 2026-05-15): hermes
# inherits TimeoutStopSec=4min from the systemd defaults and hangs in
# 'deactivating' for the full window when SIGTERM lands mid-tool-call.
# With the systemwide `TimeoutStopFailureMode=abort` drop-in in place,
# systemd then SIGABRTs the process at the 4-minute mark -- net effect:
# every restart blocks the operator's terminal for ~4 minutes when the
# gateway is mid-anything.
#
# 60 s pairs with hermes's internal agent.restart_drain_timeout=20
# (mios-hermes-firstboot seeds this in config.yaml). Hermes's
# invariant: TimeoutStopSec must be >= drain_timeout + 30 s buffer
# for graceful drain + cleanup. With drain=20 + 30 buffer = 50, so
# 60 s gives a small safety margin.
#
# Net shutdown latency:
#   * Typical (no in-flight work):     ~2-3 s
#   * Mid-tool-call (one OWUI chat):    20-30 s (drain finishes,
#                                                hermes exits clean)
#   * Worst case (drain hung):          60 s (SIGKILL via systemd)
#
# vs the upstream default 240 s (4 min) -- 4-8x faster restart
# cycle for the common dev case. Operator directive 2026-05-15:
# `mios-restart hermes` shouldn't take 4 minutes for a routine
# skill edit.
#
# Operators with a long-running synchronous task should wait for it
# to complete BEFORE restarting. For tasks that legitimately need
# minutes, raise agent.restart_drain_timeout in config.yaml AND
# raise TimeoutStopSec here in lock-step (drain + 30 s).

[Service]
TimeoutStopSec=60s
