# usr/lib/systemd/system/hermes-agent.service.d/60-workdir-root.conf
# /usr/lib/systemd/system/hermes-agent.service.d/60-workdir-root.conf
# Run hermes FROM the MiOS root (/) so it sees the root-level MiOS AI identity
# files (MiOS.md, AGENTS.md, system-prompt.md, CLAUDE.md, GEMINI.md) in its CWD,
# the agents.md convention. / is read-only under ProtectSystem=strict (reads
# only); hermes still writes its state to HERMES_HOME=/var/lib/mios/hermes.
# operator 2026-06-10: "hermes installed to root so it sees all the MiOS AI MD files".
[Service]
WorkingDirectory=/
