# usr/lib/systemd/system/nvidia-cdi-refresh.service.d/10-mios-ordering.conf
# 'MiOS' v0.2.4 drop-in for nvidia-cdi-refresh.service
# ----------------------------------------------------------------------------
# Workaround for NVIDIA/nvidia-container-toolkit#1735: v0.2.0 added
# `After=multi-user.target` to the stock unit, creating an ordering cycle
# with any unit that both WantedBy=multi-user.target and
# Requires=nvidia-cdi-refresh.service (e.g. podman.service, k3s.service).
#
# Empty `After=` clears the inherited value, then re-declare minimal sane
# ordering: kernel modules loaded + udev coldplug fired (udev-trigger, not
# the deprecated udev-settle which warns at boot).
# ----------------------------------------------------------------------------
[Unit]
After=
After=systemd-modules-load.service systemd-udev-trigger.service
