# usr/lib/systemd/system/nvidia-cdi-refresh.service.d/10-mios-ordering.conf
# MiOS v0.2.0 drop-in for nvidia-cdi-refresh.service
# ----------------------------------------------------------------------------
# Workaround for NVIDIA/nvidia-container-toolkit#1735: v0.2.0 added
# `After=multi-user.target` to the stock unit, creating an ordering cycle
# with any unit that both WantedBy=multi-user.target and
# Requires=nvidia-cdi-refresh.service (e.g. podman.service, k3s.service).
#
# Empty `After=` clears the inherited value, then we re-declare minimal
# sane ordering: wait for kernel modules loaded and udev settled.
# ----------------------------------------------------------------------------
[Unit]
After=
After=systemd-modules-load.service systemd-udev-settle.service
Wants=systemd-udev-settle.service
