# usr/lib/systemd/system/mios-freeipa-enroll.service.d/10-boot-timeout.conf
# AI-hint: Configures a 120-second startup timeout for the mios-freeipa-enroll service to prevent system hangs during FreeIPA enrollment.
# AI-related: mios-freeipa-enroll
[Service]
TimeoutStartSec=120
