# usr/lib/systemd/system/libvirtd.service.d/10-mios.conf
# AI-hint: Configures libvirtd.service to ensure it starts after the libvirtd.socket and extends the shutdown timeout to 120s to prevent race conditions with Cockpit and handle slow VM storage during teardown.
# AI-related: libvirtd.service, libvirtd.socket
# Prevent startup race condition between Cockpit and libvirt on ucore bases
[Unit]
After=libvirtd.socket

[Service]
# Increase shutdown timeout for large VMs / slow storage
TimeoutStopSec=120
