# usr/lib/systemd/system/libvirtd.service.d/10-mios.conf
# Prevent startup race condition between Cockpit and libvirt on ucore bases
[Unit]
After=libvirtd.socket

[Service]
# Increase shutdown timeout for large VMs / slow storage
TimeoutStopSec=120
