# usr/lib/systemd/system/libvirtd.service.d/override.conf
# 'MiOS': prevent SIGKILL of slow VMs during shutdown
[Service]
TimeoutStopSec=120s