# usr/lib/systemd/system/libvirtd.service.d/override.conf
# AI-hint: Overrides the libvirtd service timeout to 120s to prevent the systemd manager from sending SIGKILL to slow-running virtual machines during the system shutdown sequence.
# 'MiOS': prevent SIGKILL of slow VMs during shutdown
[Service]
TimeoutStopSec=120s