# usr/lib/systemd/system/mios-llm-heavy.service.d/50-blade-gpu-serving.conf
# AI-hint: GENERATED systemd capability drop-in for MiOS (WS-BLADE). DO NOT EDIT -- regenerate via tools/generate-blade-dropins.py.
# Gates this unit under systemd using the /etc/mios/blade.d/gpu-serving marker.
[Unit]
ConditionPathExists=/etc/mios/blade.d/gpu-serving
