# usr/lib/systemd/system/systemd-binfmt.service.d/10-wsl.conf
# AI-hint: Prevents systemd-binfmt.service from entering a failed state under WSL2 when /proc/sys/fs/binfmt_misc is mounted read-only.
[Unit]
ConditionPathIsReadWrite=/proc/sys/fs/binfmt_misc
