# usr/lib/systemd/user/mios-launcher.service
# AI-hint: Systemd user unit for the mios-launcher-daemon; acts as a privileged broker that listens on a Unix socket to execute GUI applications with the operator's environment (WSLg/Wayland) for agents lacking execution rights.
# AI-related: /usr/libexec/mios/mios-launcher-daemon, mios-launcher-daemon, mios-launcher, mios-hermes, mios-launch, mios-gui, mios-os-control, default.target
[Unit]
# /usr/lib/systemd/user/mios-launcher.service
#
# Operator-side launcher broker. Runs in user@<uid>.service so it
# inherits the operator's WSLg env (WAYLAND_DISPLAY,
# WSL2_GUI_APPS_ENABLED=1, WSL_INTEROP, DBUS_SESSION_BUS_ADDRESS,
# the live wayland socket on /run/user/<uid>/wayland-0). Listens on
# a unix socket the agent (mios-hermes) can write launch requests
# to, dispatching them as the operator with the right environment.
#
# Why this exists: mios-hermes (uid 820) cannot exec
# /mnt/c/Windows/System32/*.exe (mode 0544 owned by mios:mios via
# WSL metadata; even root chmod fails because Windows ACL blocks
# the xattr write). The broker is the operator's standing service
# that DOES have exec access; the agent posts launch requests to
# its socket from any context.
#
# Operator directive 2026-05-15: "MiOS-Agents should be able to do
# this all again!! FIX!!! JUST CREATE universal Launching tools/
# skills for launching applications anywhere in any environment(s)
# (cross-platform)".
Description='MiOS' launcher broker (operator-side)
Documentation=man:mios-launcher-daemon
#
# Run ONLY in the OPERATOR's user manager. The broker binds the shared
# /run/mios-launcher/launcher.sock; if root's (uid 0) or a second user's
# (uid 1000) user manager ALSO starts it they contend for that socket, and
# whichever wins launches GUI apps in the WRONG session (root -> invisible,
# then the window dies). ConditionUser gates the broker to the operator so
# exactly ONE instance owns the socket. ('mios' = MIOS_USER, the operator;
# render from MIOS_USER if that default ever changes.) Operator-hit
# 2026-06-06: multi-broker invisible launches recurring after WSL restarts.
ConditionUser=mios

[Service]
Type=simple
# mios-launch + the verbs it shells to (flatpak-launch, mios-gui, mios-os-control)
# live in /usr/libexec/mios, which is NOT on the default systemd-user PATH -- so the
# broker got "mios-launch: command not found" (exit 127) and EVERY app launch failed
# (operator 2026-06-05 "open epiphany"). Put libexec + the mios shim dir on PATH.
Environment=PATH=/usr/libexec/mios:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/usr/libexec/mios/mios-launcher-daemon
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
