# usr/lib/systemd/user/mios-launcher.service
[Unit]
# /usr/lib/systemd/user/mios-launcher.service
#
# Operator-side launcher broker. Runs in user@<uid>.service so it
# inherits the operator's WSLg env (WAYLAND_DISPLAY,
# WSL2_GUI_APPS_ENABLED=1, WSL_INTEROP, DBUS_SESSION_BUS_ADDRESS,
# the live wayland socket on /run/user/<uid>/wayland-0). Listens on
# a unix socket the agent (mios-hermes) can write launch requests
# to, dispatching them as the operator with the right environment.
#
# Why this exists: mios-hermes (uid 820) cannot exec
# /mnt/c/Windows/System32/*.exe (mode 0544 owned by mios:mios via
# WSL metadata; even root chmod fails because Windows ACL blocks
# the xattr write). The broker is the operator's standing service
# that DOES have exec access; the agent posts launch requests to
# its socket from any context.
#
# Operator directive 2026-05-15: "MiOS-Agents should be able to do
# this all again!! FIX!!! JUST CREATE universal Launching tools/
# skills for launching applications anywhere in any environment(s)
# (cross-platform)".
Description='MiOS' launcher broker (operator-side)
Documentation=man:mios-launcher-daemon

[Service]
Type=simple
# mios-launch + the verbs it shells to (flatpak-launch, mios-gui, mios-os-control)
# live in /usr/libexec/mios, which is NOT on the default systemd-user PATH -- so the
# broker got "mios-launch: command not found" (exit 127) and EVERY app launch failed
# (operator 2026-06-05 "open epiphany"). Put libexec + the mios shim dir on PATH.
Environment=PATH=/usr/libexec/mios:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
# WSLg display env (operator 2026-06-06: epiphany launched but no window appeared --
# the broker had NO WAYLAND_DISPLAY/XDG_RUNTIME_DIR in its env, so every GUI app it
# spawns runs headless with no compositor and exits). Point it at the live WSLg
# socket (/run/user/992/wayland-0 -> /mnt/wslg/runtime-dir/wayland-0) so launched
# apps surface on the operator's desktop. uid 1000 = the broker user.
Environment=XDG_RUNTIME_DIR=/run/user/992
Environment=WAYLAND_DISPLAY=wayland-0
Environment=DISPLAY=:0
Environment=WSL2_GUI_APPS_ENABLED=1
# flatpak run (epiphany) needs the dbus SESSION bus for the portal/session or it
# exits (exit 2, no window) -- operator 2026-06-06: broker ran mios-launch epiphany
# but no window. Point at uid 1000's session bus.
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/992/bus
# Use the DIRECT launch path (skip flatpak-launch's wsl.exe-session re-launch):
# (1) the nested wsl.exe session tears down its children after the bootstrap so
# apps die; (2) X11 :0 is broken on this WSLg (xdg-desktop-portal-gtk: "cannot
# open display :0") so launches must use WAYLAND. The broker is a persistent
# user@992 service with the wayland-0 socket + user session-bus + portal, so a
# setsid direct launch survives and reaches the wayland display (operator
# 2026-06-06). Pair with the per-app override GDK_BACKEND=wayland.
Environment=_MIOS_FLATPAK_INSESSION=1
ExecStart=/usr/libexec/mios/mios-launcher-daemon
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
