# usr/lib/systemd/user/mios-wsl-graphical-session.service
[Unit]
Description=MiOS: pull in graphical-session.target on WSLg
Documentation=https://github.com/mios-dev/MiOS
Documentation=https://systemd.io/DESKTOP_ENVIRONMENTS/

# WSLg-only. On bare-metal MiOS the GNOME session manager
# (gnome-session.service) is the canonical activator of
# graphical-session.target -- we don't want to fight that here.
# On WSLg there is NO gnome-session (WSLg's compositor is the
# Wayland server, NOT a GNOME shell), so graphical-session.target
# never activates and every BindsTo=graphical-session.target unit
# fails -- including xdg-desktop-portal-gnome -> a "Dependency
# failed" cascade to xdg-desktop-portal -> nautilus / epiphany /
# gnome-software / every flatpak that touches a portal break.
#
# Canonical fix per upstream conventions (Sway/Hyprland model,
# Arch waybar wiki, systemd.io/DESKTOP_ENVIRONMENTS): a oneshot
# user unit with Wants=graphical-session.target and
# WantedBy=default.target. Dependency-driven activation pulls in
# graphical-session.target legitimately -- RefuseManualStart=yes
# only blocks `systemctl start graphical-session.target`, not
# Wants= dependency activation.
#
# Env-import precondition: xdg-desktop-portal-gnome refuses to
# operate without XDG_CURRENT_DESKTOP / WAYLAND_DISPLAY /
# DISPLAY / XDG_SESSION_TYPE in the systemd-user environment.
# WSLg sets them in the login shell only, so we run
# mios-wsl-env-import.service first to push them into systemd
# user-bus + dbus activation environment.
ConditionPathIsDirectory=/mnt/wslg
DefaultDependencies=no
After=basic.target
After=mios-wsl-env-import.service
Requires=mios-wsl-env-import.service
Wants=graphical-session.target
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=default.target
