# etc/systemd/system/hermes-agent.service.d/50-soul-sync.conf
# /etc/systemd/system/hermes-agent.service.d/50-soul-sync.conf
#
# Re-sync /var/lib/mios/hermes/SOUL.md from
# /usr/share/mios/ai/hermes-soul.md on every service start so SOUL.md
# updates propagate without a manual touch+restart. Operator-flagged
# 2026-05-17: WSLg rule didn't take effect after my deploy because
# `cp` preserved mtime, so the firstboot-style content-diff check
# said "no change" -- the new prose never reached the agent.
#
# mios-hermes-soul-sync exits 0 on noop OR successful sync, so this
# ExecStartPre is safe to keep enabled even when SOUL.md is current.
# It also preserves any operator-appended runtime-context block (the
# ### MIOS-RUNTIME-CONTEXT-BEGIN section firstboot adds for Discord
# defaults, etc.).

[Service]
ExecStartPre=/usr/libexec/mios/mios-hermes-soul-sync
