# usr/lib/systemd/system/cockpit-wsinstance-socket-user.service.d/10-mios-container.conf
# /usr/lib/systemd/system/cockpit-wsinstance-socket-user.service.d/10-mios-container.conf
#
# Cockpit's wsinstance-socket-user service uses DynamicUser=yes which
# pulls in PrivateTmp / PrivateMounts / mount-namespacing -- those
# require kernel privileges (CAP_SYS_ADMIN + namespaces) that are not
# available inside a container or a podman-machine WSL distro. The
# stock unit then fails with:
#
#     Failed to set up mount namespacing: Operation not supported
#     Failed at step NAMESPACE spawning /bin/true: Operation not supported
#     status=226/NAMESPACE
#
# The systemd start-rate-limiter then trips on cockpit-wsinstance-
# socket-user.service after a handful of retries, which cascades:
# cockpit-wsinstance-http.socket fails -> cockpit.service fails ->
# cockpit.socket trigger-limit-hit -> the operator's dashboard shows
# 'Cockpit https://localhost:9090/  (open circle = inactive)'.
#
# This drop-in is gated to ConditionVirtualization=container so it
# applies inside MiOS-DEV (podman-machine) but a NO-OP on bare-metal
# / Hyper-V / QEMU MiOS deployments where the namespacing works
# correctly.

[Unit]
ConditionVirtualization=container

[Service]
# Strip every namespace-requiring directive. DynamicUser=no + a static
# UnsetEnvironment block let the unit actually start; the /run/cockpit/
# wsinstance/ socket dir falls back to the cockpit user that the RPM
# already creates statically.
DynamicUser=no
PrivateTmp=no
PrivateMounts=no
PrivateDevices=no
PrivateUsers=no
ProtectHome=no
ProtectSystem=no
ProtectControlGroups=no
ProtectKernelModules=no
ProtectKernelTunables=no
ProtectKernelLogs=no
RestrictNamespaces=no
RestrictRealtime=no
LockPersonality=no
NoNewPrivileges=no
SystemCallFilter=
User=cockpit-ws
Group=cockpit-ws
