# etc/profile.d/gnupg2.sh
if [[ $- == *i* ]]; then
    export GPG_TTY=$(tty)
fi
