# etc/profile.d/gnupg2.csh
if ( $?prompt ) then
    setenv GPG_TTY `tty`
endif
