# usr/lib/systemd/system/cockpit.socket.d/10-mios-wsl2.conf
[Socket]
# WSL2: cockpit.socket works fine under systemd-in-WSL2.
# No ConditionVirtualization block — allow socket to start.
FreeBind=yes
