# usr/lib/systemd/system/cockpit.socket.d/10-mios-wsl2.conf
[Unit]
# In WSL2, remove the libvirtd.socket ordering -- libvirtd works fine
# but the ordering causes an unnecessary delay on first connection.
After=
After=dbus-daemon-wsl.service
