# usr/lib/systemd/system/cockpit.socket.d/10-mios.conf
# Prevent cockpit.socket from activating before libvirtd.socket is ready.
# Without this, opening the Machines panel immediately after boot races
# against libvirtd startup and produces "Failed to connect to libvirt" errors.
[Unit]
After=libvirtd.socket
