# usr/lib/systemd/system/cockpit.socket.d/10-mios-wsl2.conf
# AI-hint: systemd drop-in for cockpit.socket adding FreeBind=yes so the Cockpit web console (:9090) can bind under WSL2 even before the address is up; intentionally sets no ConditionVirtualization so the socket starts in WSL.
# AI-related: cockpit.socket
[Socket]
# WSL2: cockpit.socket works fine under systemd-in-WSL2.
# No ConditionVirtualization block -- allow socket to start.
FreeBind=yes
