# usr/lib/systemd/system/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
# Explicit 0.0.0.0:9090 (IPv4) -- the bare `ListenStream=9090` form
# resolves to [::]:9090 dual-stack on dual-stack-by-default Linux,
# which appears as `*:9090` in ss output. WSL2's NAT-mode
# localhostForwarding does NOT forward dual-stack binds; only
# explicit IPv4 listens reach Windows-side 127.0.0.1:9090. Operator-
# confirmed 2026-05-11: cockpit was active inside MiOS-DEV but
# Zen browser timed out on https://localhost:9090 from Windows
# until this drop-in bound 0.0.0.0:9090 explicitly.
ListenStream=0.0.0.0:9090
FreeBind=yes
# WSL2 Windows-side port forwarder (wslhost.exe) keeps 9090 mapped for
# a short window after the Linux process exits, so the next boot of
# the same distro hits EADDRINUSE on cockpit.socket and the unit fails
# with "Failed with result 'resources'". ReusePort lets the kernel
# bind anyway -- safe because cockpit.socket is the only listener
# that ever wants 9090 inside the distro.
ReusePort=yes
