# usr/lib/systemd/system/cockpit.socket.d/listen.conf
# AI-hint: Configures the cockpit.socket drop-in to bind specifically to 0.0.0.0:9090 and enable ReusePort, ensuring IPv4 accessibility from Windows hosts via WSL2 and preventing EADDRINUSE errors during rapid restarts.
# AI-related: cockpit.socket, localhost:9090
[Socket]
ListenStream=
# Explicit 0.0.0.0:<port> (IPv4) -- the bare `ListenStream=<port>` form
# resolves to [::]:<port> dual-stack on dual-stack-by-default Linux,
# which appears as `*:<port>` in ss output. WSL2's NAT-mode
# localhostForwarding does NOT forward dual-stack binds; only explicit
# IPv4 listens reach the Windows-side 127.0.0.1:<port>.
# systemd does NOT expand ${VAR} / ${VAR:-default} in ListenStream -- it
# fails "Failed to parse address ... Invalid argument" and the socket
# then refuses with "Unit has no Listen setting", so cockpit never binds.
# The port is therefore a LITERAL mirroring mios.toml [ports].cockpit
# (MIOS_PORT_COCKPIT); keep the two in sync.
ListenStream=0.0.0.0:8090
FreeBind=yes
# WSL2 Windows-side port forwarder (wslhost.exe) keeps 9090 mapped for
# a short window after the Linux process exits, so the next boot of
# the same distro hits EADDRINUSE on cockpit.socket and the unit fails
# with "Failed with result 'resources'". ReusePort lets the kernel
# bind anyway -- safe because cockpit.socket is the only listener
# that ever wants 9090 inside the distro.
ReusePort=yes
