# usr/lib/systemd/system/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=9090
FreeBind=yes
# WSL2 Windows-side port forwarder (wslhost.exe) keeps 9090 mapped for
# a short window after the Linux process exits, so the next boot of
# the same distro hits EADDRINUSE on cockpit.socket and the unit fails
# with "Failed with result 'resources'". ReusePort lets the kernel
# bind anyway -- safe because cockpit.socket is the only listener
# that ever wants 9090 inside the distro.
ReusePort=yes
