# usr/lib/systemd/system/waydroid-container.service.d/10-nvidia-swiftshader.conf
[Service]
# Dynamically configure EGL software rendering for NVIDIA hosts before container start
ExecStartPre=/usr/libexec/mios/mios-waydroid-fallback.sh
