# usr/lib/systemd/system/hermes-agent.service
[Unit]
# /usr/lib/systemd/system/hermes-agent.service
#
# Hermes-Agent as a DIRECT host service -- not a container Quadlet.
# Operator directive 2026-05-13: the agent is installed natively onto
# the MiOS root filesystem (automation/38-hermes-agent.sh) and runs
# here as a host systemd service so it sees the real `/`, the real
# systemd, and the real git working tree without container bind-mount
# indirection. The mios-hermes + mios-hermes-dashboard CONTAINER
# Quadlets are shelved (mios.toml [quadlets.enable] = false).
#
# Inference backend: the OpenAI-compatible /v1 surface resolved from
# mios.toml [ai].endpoint (default http://localhost:11434/v1 = the
# ollama container; swap to vllm / llama.cpp by editing mios.toml).
# This unit deliberately depends on ollama.service so the backend is
# up first, but stays Wants= (not Requires=) so the agent still starts
# -- and can be repointed -- if ollama is swapped out.
Description='MiOS' Hermes-Agent (direct host install, OpenAI gateway)
Documentation=https://github.com/NousResearch/hermes-agent
After=network-online.target ollama.service
Wants=network-online.target ollama.service
# Skip cleanly if the direct install didn't land (e.g. clone failed at
# build time) rather than crash-looping.
ConditionPathExists=/usr/lib/mios/hermes-agent/.venv/bin/hermes

[Service]
Type=simple
# Runtime identity: the mios-hermes sysuser (uid 820) already exists in
# /usr/lib/sysusers.d/ from the container era -- reuse it so /var/lib
# ownership stays consistent across the container->direct migration.
User=mios-hermes
Group=mios-hermes
# Runtime state dir -- systemd creates + chowns it (NO-MKDIR-IN-VAR law).
# ConfigurationDirectory= deliberately omitted: /etc/mios/hermes is
# pre-created by mios-hermes-firstboot at mode 0750; declaring
# ConfigurationDirectory=mios/hermes here makes systemd fight that
# (warns "mode is different. File system: 750 ConfigurationDirectoryMode:
# 755"). ReadWritePaths=/etc/mios below already grants write access.
StateDirectory=mios/hermes
Environment=HOME=/var/lib/mios/hermes
Environment=HERMES_HOME=/var/lib/mios/hermes
Environment=HERMES_CONFIG_PATH=/etc/mios/hermes/config.yaml
# Gateway bind: PORT is the OpenAI /v1 surface port (mios.toml
# [ports].hermes, default 8642). API_SERVER_HOST=0.0.0.0 comes from
# api.env (so the WSL2 portproxy + LAN reach it). HERMES_BACKEND_BASE_URL
# points at the ollama container's host-net OpenAI surface -- hermes-agent
# runs ON the host now (direct install), so localhost:11434 reaches the
# Network=host ollama container directly.
Environment=PORT=8642
Environment=HERMES_BACKEND_BASE_URL=http://localhost:11434
# api.env carries API_SERVER_HOST=0.0.0.0, API_SERVER_KEY, the dashboard
# token, GATEWAY_ALLOW_ALL_USERS, etc. install.env is the non-secret
# bootstrap bridge. Both optional (- prefix) so the unit still starts
# pre-firstboot.
EnvironmentFile=-/etc/mios/install.env
EnvironmentFile=-/etc/mios/hermes/api.env
# `hermes gateway run` serves the OpenAI /v1 surface + agent dashboard
# in the foreground (Type=simple). --replace evicts any stale gateway
# instance -- essential for clean systemd restarts. Host/port come from
# the env above (API_SERVER_HOST + PORT), NOT CLI flags: `hermes gateway`
# itself is a service-manager subcommand whose only positional is
# {run,start,stop,...} -- passing --host/--port to it is an arg error.
ExecStart=/usr/lib/mios/hermes-agent/.venv/bin/hermes gateway run --replace
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
# hermes gateway's restart drain (agent.restart_drain_timeout, default
# 180s) needs TimeoutStopSec >= drain_timeout + 30s headroom or systemd
# SIGKILLs the gateway mid-drain (operator-confirmed 2026-05-14 warning:
# "TimeoutStopSec=45s but drain_timeout=180s, expected >=210s").
TimeoutStopSec=240s
# Hardening -- the agent reads `/` but should not be able to rewrite
# vendor /usr; /etc + /var/lib/mios/hermes stay writable via the
# Directory= directives above.
ProtectSystem=strict
ReadWritePaths=/etc/mios /var/lib/mios
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
