# usr/lib/systemd/system/hermes-agent.service
[Unit]
# /usr/lib/systemd/system/hermes-agent.service
#
# Hermes-Agent as a DIRECT host service -- not a container Quadlet.
# Operator directive 2026-05-13: the agent is installed natively onto
# the MiOS root filesystem (automation/38-hermes-agent.sh) and runs
# here as a host systemd service so it sees the real `/`, the real
# systemd, and the real git working tree without container bind-mount
# indirection. The mios-hermes + mios-hermes-dashboard CONTAINER
# Quadlets are shelved (mios.toml [quadlets.enable] = false).
#
# Inference backend: the OpenAI-compatible /v1 surface resolved from
# mios.toml [ai].endpoint (default http://localhost:11434/v1 = the
# ollama container; swap to vllm / llama.cpp by editing mios.toml).
# This unit deliberately depends on ollama.service so the backend is
# up first, but stays Wants= (not Requires=) so the agent still starts
# -- and can be repointed -- if ollama is swapped out.
Description='MiOS' Hermes-Agent (direct host install, OpenAI gateway)
Documentation=https://github.com/NousResearch/hermes-agent
After=network-online.target ollama.service
Wants=network-online.target ollama.service
# Skip cleanly if the direct install didn't land (e.g. clone failed at
# build time) rather than crash-looping.
ConditionPathExists=/usr/lib/mios/hermes-agent/.venv/bin/hermes

[Service]
Type=simple
# Runtime identity: the mios-hermes sysuser (uid 820) already exists in
# /usr/lib/sysusers.d/ from the container era -- reuse it so /var/lib
# ownership stays consistent across the container->direct migration.
User=mios-ai
Group=mios-ai
# Runtime state dir: /var/lib/mios/hermes. StateDirectory= is
# deliberately NOT used -- systemd re-applies StateDirectoryMode (and
# silently drops the setgid bit, landing 0700) on EVERY start, which
# races and clobbers the 2770 group-shared mode that mios-hermes-
# firstboot sets so the `hermes` CLI (run as the mios operator, a
# mios-hermes group member) and this gateway service can share one
# config + state tree. mios-hermes-firstboot.service is ordered
# Before= this unit and does the mkdir + chown 820:820 + chmod 2770,
# so the directory is guaranteed present and correctly-permissioned
# before ExecStart. ReadWritePaths=/var/lib/mios below grants write.
# ConfigurationDirectory= likewise omitted: /etc/mios/hermes is
# pre-created by firstboot at mode 0750.
Environment=HOME=/var/lib/mios/hermes
Environment=HERMES_HOME=/var/lib/mios/hermes
# HERMES_CONFIG_PATH deliberately NOT set: Hermes-Agent 0.13.x ignores
# it entirely and reads $HERMES_HOME/config.yaml. Keeping a stale
# HERMES_CONFIG_PATH=/etc/mios/hermes/config.yaml here just misleads
# operators into editing a file the agent never loads.
# SearXNG endpoint for the agent's web-search/extract tools. Host-net
# dev-VM value; bridge builds get mios-searxng:8080 via a Quadlet-style
# drop-in or 15-render. Web search is retrieval, not inference (Law 5).
Environment=SEARXNG_URL=http://localhost:8888
# crawl-engine endpoint for the agent's `crawl` tool (mios-crawl POSTs here).
# The crawl engine is the loopback venv service mios-crawl4ai.service
# (crawl4ai over Chrome CDP + Camoufox fallback); same retrieval-not-inference
# posture as SearXNG (a local browser is not a vendor cloud model -- Law 5).
# mios-crawl defaults to this same URL, so this line is for parity +
# discoverability. (Renamed 2026-05-24: MIOS_CRAWL4AI_URL -> MIOS_CRAWL_SERVICE_URL
# when the container was scrapped for the venv service.)
Environment=MIOS_CRAWL_SERVICE_URL=http://127.0.0.1:11235
# Gateway bind: PORT is the OpenAI /v1 surface port (mios.toml
# [ports].hermes, default 8642). API_SERVER_HOST=0.0.0.0 comes from
# api.env (so the WSL2 portproxy + LAN reach it). HERMES_BACKEND_BASE_URL
# points at the ollama container's host-net OpenAI surface -- hermes-agent
# runs ON the host now (direct install), so localhost:11434 reaches the
# Network=host ollama container directly.
Environment=PORT=8642
Environment=HERMES_BACKEND_BASE_URL=http://localhost:11434
# api.env carries API_SERVER_HOST=0.0.0.0, API_SERVER_KEY, the dashboard
# token, GATEWAY_ALLOW_ALL_USERS, etc. install.env is the non-secret
# bootstrap bridge. Both optional (- prefix) so the unit still starts
# pre-firstboot.
EnvironmentFile=-/etc/mios/install.env
EnvironmentFile=-/etc/mios/hermes/api.env
# Discord bot token (DISCORD_BOT_TOKEN) for the discord toolset. The
# file is operator-managed; if absent or empty, discord_tool surfaces
# a clear "missing token" error to the agent on first call -- which
# the agent then asks the operator for. Drop the token in like:
#   echo 'DISCORD_BOT_TOKEN=...' | sudo install -m 0640 -o root \
#       -g mios-hermes /dev/stdin /etc/mios/hermes/discord.env
EnvironmentFile=-/etc/mios/hermes/discord.env
# `hermes gateway run` serves the OpenAI /v1 surface + agent dashboard
# in the foreground (Type=simple). --replace evicts any stale gateway
# instance -- essential for clean systemd restarts. Host/port come from
# the env above (API_SERVER_HOST + PORT), NOT CLI flags: `hermes gateway`
# itself is a service-manager subcommand whose only positional is
# {run,start,stop,...} -- passing --host/--port to it is an arg error.
# ExecStartPre: idempotently apply the MiOS Discord progressive-
# reactions patch to the gateway's site-packages discord.py. The
# script is a no-op when the marker comment is already present, so
# every restart re-asserts the patch -- survives venv refreshes
# (pip install -e .) that would otherwise drop the modification.
# Exit codes: 0 = applied or already-applied; 1 = file missing
# (skipped, service still starts); 2 = upstream refactored (warn,
# service still starts). The `-` prefix makes systemd ignore the
# exit code so a Discord-gateway-only failure can't block the
# whole agent from booting.
# Guarantee the runtime state dir is OWNED BY the service user (mios-ai)
# on EVERY start. The decisive fix is the ownership: the `hermes gateway`
# tightens HERMES_HOME to 0700 itself on startup, so as long as mios-ai
# owns the dir it can create gateway.pid and manage its own tree. The `+`
# prefix runs this as root (User=mios-ai cannot chown); `install -d` is
# idempotent -- on the existing dir it just re-asserts owner/group/mode.
# 2026-05-24 incident: /var/lib/mios/hermes had drifted to root-era
# ownership (mios-hermes:mios-hermes, mode 0755) AFTER firstboot (which
# exits early on later boots once config exists, so it no longer
# re-corrects the dir). The mios-ai gateway, NOT the owner, could not
# create gateway.pid -> "PermissionError: [Errno 13]" which the gateway
# misreports as "PID file race lost to another gateway instance" -> a
# crash/restart loop with NOTHING listening on :8642 (silently "active
# (running)" but a dead agent). 0700 here matches what the gateway sets
# anyway; owner mios-ai retains full access.
ExecStartPre=+/usr/bin/install -d -o mios-ai -g mios-ai -m 0700 /var/lib/mios/hermes
ExecStartPre=-/usr/bin/python3 /usr/libexec/mios/mios-hermes-discord-reactions-patch /usr/lib/mios/hermes-agent/.venv/lib/python3.14/site-packages/gateway/platforms/discord.py
ExecStart=/usr/lib/mios/hermes-agent/.venv/bin/hermes gateway run --replace
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s
# hermes gateway's restart drain (agent.restart_drain_timeout, default
# 180s) needs TimeoutStopSec >= drain_timeout + 30s headroom or systemd
# SIGKILLs the gateway mid-drain (operator-confirmed 2026-05-14 warning:
# "TimeoutStopSec=45s but drain_timeout=180s, expected >=210s").
TimeoutStopSec=240s
# File-descriptor ceiling. systemd's default soft RLIMIT_NOFILE for a
# service is 1024 -- far too low for a long-running async gateway that
# holds many concurrent sockets at once (ollama backend, Discord
# gateway, the agent dashboard, and the swarm's every-5s health probes).
# 2026-05-24 incident: the gateway accumulated fds over a 15h uptime
# until it hit 1024, after which asyncio's accept() failed with EMFILE
# ("OSError: [Errno 24] Too many open files" / "socket.accept() out of
# system resource"). The unit stayed "active (running)" but could no
# longer accept ANY connection on :8642 (95 conns jammed in the listen
# backlog) so every chat/tool call timed out -- a silent, wedged agent.
# Raise the soft limit to the system hard cap (1048576) so a slow fd
# leak can no longer wedge the agent; a restart still resets the count.
LimitNOFILE=1048576
# Hardening -- the agent reads `/` but should not be able to rewrite
# vendor /usr; /etc + /var/lib/mios/hermes stay writable via the
# Directory= directives above.
ProtectSystem=strict
ReadWritePaths=/etc/mios /var/lib/mios /run/mios-launcher
ProtectHome=false
# NoNewPrivileges=false (was true) -- LOAD-BEARING for the operator's
# explicit sudo authorization. With NoNewPrivileges=true the kernel
# refuses to honour sudo's SUID bit (or any setuid binary) regardless
# of /etc/sudoers.d/, so `sudo mios build` from the agent dies with
# "sudo: The "no new privileges" flag is set, which prevents sudo from
# running as root" -- the sudoers grant becomes a no-op. The operator
# explicitly authorized mios-hermes to sudo (see /etc/sudoers.d/
# 10-mios-hermes); keeping NoNewPrivileges=true would silently
# undermine that authorization. Revoke either by deleting the sudoers
# drop-in OR by flipping this back to true; both correctly trip the
# build driver's privilege preflight (operator-confirmed 2026-05-15).
NoNewPrivileges=false
# PrivateTmp=yes is LOAD-BEARING. ProtectSystem=strict makes /tmp and
# /var/tmp READ-ONLY for this service. Hermes-Agent's `terminal` tool
# (and Python's tempfile underneath every code_execution / shell call)
# needs a writable temp dir to stage execution -- with none, EVERY tool
# call fails with "[Errno 2] No usable temporary directory found in
# ['/tmp', '/var/tmp', '/usr/tmp', '/']". The agent, unable to actually
# run anything, then MAKES UP plausible-looking output (fake build
# logs, invented config flags) instead of erroring honestly --
# operator-confirmed 2026-05-14, the "confused local agent" report.
# PrivateTmp=yes gives the service its own private, writable /tmp +
# /var/tmp namespace, which is exactly what tempfile needs and is more
# secure than adding the shared /tmp to ReadWritePaths. The agent's
# real work still targets /, /etc, /var/lib/mios -- /tmp is only scratch.
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
