# usr/lib/systemd/system/mios-opencode-gateway.service
[Unit]
Description='MiOS' OpenCode /v1 gateway (OpenAI adapter fronting the opencode CLI)
Documentation=file:///usr/lib/mios/opencode-gateway/server.py
After=network-online.target ollama.service
Wants=network-online.target
# Makes [agents.opencode] endpoint :8633/v1 a REAL OpenAI endpoint so
# agent-pipe's multi-agent fan-out (opencode secondary) + the primary path
# (refine target_agent=opencode) can reach it. opencode has no native /v1
# (its `serve` is opencode's own OpenAPI on :4096); this shim wraps
# `opencode run`. FOSS + offline (opencode -> local ollama qwen2.5-coder).
#
# NOTE: opencode's model runs on the dGPU lane, so under VRAM pressure a
# run times out (gateway returns empty -> fan-out drops opencode). Deploy
# when the GPU has headroom; until then the unit can be left disabled.

[Service]
Type=simple
# v1 runs as root: root holds the /usr/local/bin/opencode symlink AND the
# local-ollama provider config (/root/.config/opencode/opencode.json).
# Hardening to a dedicated mios-opencode sysuser + ProtectSystem=strict is
# a follow-up (matches the mios-agent-pipe.service hardening pattern).
Environment=MIOS_PORT_OPENCODE_GATEWAY=8633
Environment=MIOS_OPENCODE_BIN=/usr/local/bin/opencode
Environment=MIOS_OPENCODE_WORKDIR=/var/lib/mios/opencode-gateway/work
Environment=MIOS_OPENCODE_TIMEOUT_S=90
Environment=HOME=/root
# Reuse the hermes-agent venv interpreter (fastapi/uvicorn/starlette already
# installed there) -- same pattern as mios-agent-pipe.service.
ExecStart=/usr/lib/mios/hermes-agent/.venv/bin/python3 -u /usr/lib/mios/opencode-gateway/server.py
Restart=on-failure
RestartSec=10s
TimeoutStartSec=30s

[Install]
WantedBy=multi-user.target
