# usr/lib/systemd/system/mios-podman-gc.timer
[Unit]
Description=Weekly Podman Cleanup
[Timer]
OnCalendar=weekly
Persistent=true
[Install]
WantedBy=timers.target
