# usr/lib/systemd/system/mios-agent-pipe.service
# AI-hint: Systemd unit defining the agent-pipe FastAPI service which acts as a router, refiner, and critic for the Hermes gateway, routing inference requests to the llama.cpp backend on port 11450.
# AI-related: /usr/lib/mios/agent-pipe/server.py, /etc/mios/hermes/api.env, /etc/mios/agent-pipe.env, /etc/mios/agent-pipe.env., /usr/lib/mios/agents/.venv/bin/python3, mios-surrealdb, mios-passport-provision, mios-ai, mios-hermes, mios-env
[Unit]
Description='MiOS' Agent Pipe (router + refine + critic FastAPI; fronts hermes for every gateway)
Documentation=file:///usr/lib/mios/agent-pipe/server.py
After=network-online.target hermes-agent.service mios-pgvector.service mios-passport-provision.service
Wants=network-online.target hermes-agent.service mios-passport-provision.service

[Service]
Type=simple
User=mios-ai
Group=mios-ai
# Run with the mios-hermes supplementary group so the service can read
# /etc/mios/hermes/api.env (mode 0640, group mios-hermes) and load the
# backend API_SERVER_KEY. Without it, _BACKEND_KEY stays empty and KEYLESS
# direct callers (Firefox Smart Window, curl, MCP) reach Hermes with no auth
# and get 401 -- only the OWUI path worked because it forwards the operator
# token. The agent-pipe fronts Hermes, so sharing its key group is intended.
SupplementaryGroups=mios-hermes
WorkingDirectory=/var/lib/mios/agent-pipe

# SSOT chain: every operator-tunable knob sources from mios.toml via
# /etc/profile.d/mios-env.sh, which exports MIOS_* into the service
# environment. The service code reads MIOS_PORT_AGENT_PIPE +
# MIOS_AGENT_PIPE_BACKEND + MIOS_DB_URL etc. as os.environ overrides.
EnvironmentFile=-/etc/mios/agent-pipe.env
Environment=MIOS_PORT_AGENT_PIPE=8640
# WS-10 (operator 2026-06-05 "EVERYTHING IS LLAMA.CPP, NO OLLAMA"): the agent-pipe
# reasons DIRECTLY on mios-llm-light (:11450) instead of via Hermes->ollama. All
# inference is on the dGPU via llama.cpp; ollama bypassed.
# 2026-06-15 (Claude): model fleet swapped gemma4:12b -> granite4.1:8b (the served
# brain on mios-llm-light; see mios-llm-light.yaml + mios.toml:4061). gemma4:12b is
# NO LONGER SERVED -> any role still naming it 404s on llama-swap ("no model id").
# Align EVERY pipeline role to the ONE resident served brain (granite4.1:8b): set the
# BASE default MIOS_STACK_MODEL (which _MICRO/DCI/SWARM/refine/polish/planner inherit
# via server.py) + the explicit BACKEND/SWARM/MICRO names. ONE resident model = no
# llama-swap thrash. Override per-deployment in /etc/mios/agent-pipe.env.
Environment=MIOS_AGENT_PIPE_BACKEND=http://localhost:11450/v1
Environment=MIOS_AGENT_PIPE_BACKEND_MODEL=granite4.1:8b
Environment=MIOS_STACK_MODEL=granite4.1:8b
Environment=MIOS_MICRO_MODEL=granite4.1:8b
Environment=MIOS_MICRO_ENDPOINT=http://localhost:11450/v1
Environment=MIOS_SWARM_MODEL=granite4.1:8b
# llama.cpp/mios-llm-light (:11450) is OpenAI /v1 ONLY (no ollama /api/chat, /api/ps)
# and 400s on tool_choice=required. Register it so the pipe treats it as llamacpp:
# skip forced tool_choice (no 400) + use the /slots KV-paging path. No bare port
# literal in routing code -- these are the env-SSOT hint lists.
# llama.cpp b9519 (the dGPU mios-llm-light) ACCEPTS tool_choice=required (verified) --
# only the older iGPU b9305 (:11436) rejects it. So :11450 must NOT be in the
# no-tool-choice list, or the pipe drops/auto's tool_choice and gemma4 NARRATES
# the call instead of executing it. Keep :11450 in KV_PAGING (it does /slots).
Environment=MIOS_NO_TOOL_CHOICE_HINTS=11436
Environment=MIOS_KV_PAGING_HINTS=11436,11450
Environment=MIOS_OLLAMA_API_HINTS=
# Router calls qwen3:1.7b on the CPU light-lane (:11435), not
# the dGPU/CUDA lane (:11434). Keeps router latency sub-second under
# dGPU load. Override via /etc/mios/agent-pipe.env if a deployment
# doesn't have mios-ollama-cpu standing up (e.g., bare-metal with
# only NVIDIA dGPU).
Environment=MIOS_AGENT_PIPE_ROUTER_ENDPOINT=http://localhost:11450
# CPU light-lane -> llama.cpp (was :11435 ollama, now retired). Agents that
# prefer_cpu-offload resolve here; point it at mios-llm-light so they land on
# gemma4@:11450 instead of connection-erroring on the dead ollama-cpu lane.
Environment=MIOS_OLLAMA_CPU_ENDPOINT=http://localhost:11450
# Phase A.1 -- planner runs on the dGPU/CUDA lane :11434. Falls back to
# whatever lane is reachable if :11434 is down.
# Model: mios-agent -- the canonical "MiOS AI" model (Modelfile FROM
# qwen3.5:4b as of 2026-05-22 consolidation; was gemma4:e4b). ONE light
# brain now serves plan + refine + polish with the full MiOS agent SYSTEM
# (act-don't-narrate / right-tool-for-goal / SSOT mios-os-control surface,
# ENGLISH-default, no hardcodes). The planner emits a single structured
# verb/recipe DAG; if it ever empties, DAG mode degrades to the backend
# proxy. Override via /etc/mios/agent-pipe.env.
# History: qwen2.5-coder:7b -> mios-planner(gemma4) -> mios-agent(gemma4)
# -> mios-agent(qwen3.5:4b, consolidated brain) 2026-05-22.
Environment=MIOS_AGENT_PIPE_PLANNER_ENDPOINT=http://localhost:11450
Environment=MIOS_AGENT_PIPE_PLANNER_MODEL=qwen3:1.7b
Environment=MIOS_DB_URL=http://localhost:8000
Environment=MIOS_DB_USER=root
Environment=MIOS_DB_PASS=root
Environment=MIOS_DB_NS=mios
Environment=MIOS_DB_DB=mios
# Phase C.3 -- this service signs SurrealDB writes as agent-pipe.
Environment=MIOS_PASSPORT_AGENT_NAME=agent-pipe
# Phase D.5 consolidation (operator 2026-05-22): refine (input prompt-
# enhancement) + polish (final-answer shaping WITH persona) share the ONE
# light MiOS brain -- mios-agent (FROM qwen3.5:4b, same base blob as
# mios-hermes -> no extra resident weights) on the dGPU lane :11434. The
# "polish mechanism" thus runs on a MiOS role model carrying the MiOS
# operating principles + ENGLISH-default language, not a bare base model.
# Mirrors mios.toml [refine].model / [polish].model (same pin pattern as
# MIOS_AGENT_PIPE_PLANNER_MODEL above). Override via /etc/mios/agent-pipe.env.
# REFINE runs on the LARGE 9B dGPU model (mios-hermes), POLISH on the fast 4B
# (operator 2026-05-29: "use the large dGPU model to re-write the user's queries
# with proper reasoning, intent, outcome, plan, hints/tools/skills, then that
# detailed query drives the rest of the pipeline"). Refine is THE quality lever
# -- it produces the structured plan every downstream node consumes, so it gets
# the powerful model (already dGPU-resident + warm => no extra load, just a few
# more seconds for a far better plan). Polish is output-shaping only -> stays on
# the fast 4B mios-agent so the final render is snappy. Both co-resident on the
# 4090 (9.9 + 6.5 GB). The CONCURRENT swarm secondaries still fire on OTHER
# endpoints (CPU :11435 daemon-agent, iGPU :11436, phone) at the same dispatch
# step -- the bigger refine model does not serialise them
# ([[mios_heterogeneous_parallel_dispatch]]). Override via /etc/mios/agent-pipe.env.
Environment=MIOS_REFINE_MODEL=qwen3:1.7b
Environment=MIOS_REFINE_ENDPOINT=http://localhost:11450
Environment=MIOS_POLISH_MODEL=qwen3:1.7b
Environment=MIOS_POLISH_ENDPOINT=http://localhost:11450
# Phase D.5 -- refine + polish timeouts on CPU-bound dev hardware.
# Real ROCm/CUDA finishes the LLM call in <1s; WSL CPU takes 15-30s
# for the same prompt. Generous timeout works on both.
Environment=MIOS_REFINE_TIMEOUT_S=45
Environment=MIOS_POLISH_TIMEOUT_S=45
Environment=MIOS_REFINE_MAX_TOKENS=1200

# Reuse hermes-agent's venv interpreter for fastapi/uvicorn/httpx/
# starlette -- they're already installed there. Cleaner option for a
# bootc deploy is `dnf install python3-fastapi python3-uvicorn` +
# /usr/bin/python3 (operator-confirmed 2026-05-18: python3-fastapi is
# in updates/, python3-httpx already installed). Switching to dnf
# packages is the v1 follow-up; for the scaffold commit we ride the
# hermes-agent venv to skip a package-list bump in the same change.
ExecStart=/usr/lib/mios/agents/.venv/bin/python3 -u /usr/lib/mios/agent-pipe/server.py

Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s

# Hardening: no-new-privs, private tmp, read-only system surfaces.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/mios/agent-pipe /var/lib/mios/skills /var/log /var/lib/mios/agent-env /var/lib/mios/ai
# Read-only mount of the agent-passport tree so this service can
# load its own private key + every other agent's public key for
# /passport/verify.
ReadOnlyPaths=/var/lib/mios/agent-passports

[Install]
WantedBy=multi-user.target
