# usr/lib/systemd/system/mios-flatpak-install.service
# AI-hint: Systemd unit that triggers the mios-flatpak-install script during first boot to install required Flatpak packages before the display manager starts.
# AI-related: mios-flatpak-install, display-manager.service, network-online.target, multi-user.target
[Unit]
Description='MiOS' Flatpak First-Boot Installer
ConditionVirtualization=!container
After=network-online.target
Wants=network-online.target
Before=display-manager.service
[Service]
Type=oneshot
ExecStart=/usr/libexec/mios-flatpak-install
TimeoutStartSec=600
StandardOutput=journal+console
StandardError=journal+console
[Install]
WantedBy=multi-user.target