# usr/lib/systemd/system/mios-firstboot.target
[Unit]
Description=MiOS first-boot provisioning
Documentation=https://github.com/MiOS-DEV/MiOS
After=multi-user.target
Wants=mios-cdi-detect.service mios-libvirtd-setup.service mios-grd-setup.service

[Install]
WantedBy=multi-user.target