# usr/lib/systemd/system/mios-cockpit-link.socket
# AI-hint: Socket file for mios-cockpit-link proxy. Maps port 8091 on host/WSL to the cockpit service.
# AI-related: usr/lib/systemd/system/mios-cockpit-link.service
[Unit]
Description='MiOS' Cockpit Link Proxy Socket
ConditionVirtualization=!container

[Socket]
ListenStream=0.0.0.0:8091
BindIPv6Only=both

[Install]
WantedBy=sockets.target
