# usr/lib/systemd/system/mios-libvirtd-setup.service
# AI-hint: Initializes libvirtd configuration and creates a completion flag on first boot if not in a container, ensuring virtualization services are wired and ready for MiOS.
# AI-related: /usr/libexec/mios/libvirtd-firstboot, mios-firstboot, libvirtd.socket, virtqemud.socket, virtnetworkd.socket, virtstoraged.socket, mios-firstboot.target
[Unit]
Description='MiOS' first-boot libvirtd wiring
ConditionPathExists=!/var/lib/mios/libvirtd-setup.done
ConditionVirtualization=!container
After=libvirtd.socket virtqemud.socket virtnetworkd.socket virtstoraged.socket
Wants=libvirtd.socket virtqemud.socket virtnetworkd.socket virtstoraged.socket

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mios/libvirtd-firstboot
ExecStartPost=/usr/bin/install -Dm0644 /dev/null /var/lib/mios/libvirtd-setup.done

[Install]
WantedBy=mios-firstboot.target