# usr/lib/systemd/system/mios-opencode-gateway.service
# AI-hint: Systemd unit that hosts the opencode-gateway on port 8633, providing a standard OpenAI-compatible /v1 API shim for the opencode CLI to enable multi-agent fan-out and local inference via ollama.
# AI-related: /usr/lib/mios/agents/opencode-gateway/server.py, /usr/lib/mios/agents/opencode/bin/opencode, /etc/mios/opencode/opencode.json, /usr/lib/mios/agents/.venv, /usr/lib/mios/agents/.venv/bin/python3, mios-ai, mios-agent-pipe, mios-opencode-gateway, mios-opencode, mios-llm-light.service
[Unit]
Description='MiOS' OpenCode /v1 gateway (OpenAI adapter fronting the opencode CLI)
Documentation=file:///usr/lib/mios/agents/opencode-gateway/server.py
After=network-online.target mios-llm-light.service
Wants=network-online.target mios-llm-light.service
# Makes [agents.opencode] endpoint :8633/v1 a REAL OpenAI endpoint so
# agent-pipe's multi-agent fan-out (opencode secondary) + the primary path
# (refine target_agent=opencode) can reach it. opencode has no native /v1
# (its `serve` is opencode's own OpenAPI on :4096); this shim wraps
# `opencode run`. FOSS + offline (opencode -> local ollama qwen2.5-coder).
#
# FRONT-DOOR commitment (operator 2026-05-31): opencode is a first-class /v1
# council peer via this gateway (NOT a Hermes ACP subprocess). The shim now
# passes the FULL conversation (system + history) to `opencode run`, lands its
# config at MIOS_OPENCODE_CONFIG, unifies the model id, and streams proper SSE
# deltas. Deployed + enabled by the UNIFIED agent-plane install driver
# automation/38-hermes-agent.sh: PHASE 1 builds the shared venv + Hermes,
# PHASE 2 fetches the opencode binary + lands opencode.json. The gateway
# server.py itself ships via the system-files overlay (this file's tree).
#
# This unit lives in /usr/lib/systemd/system, which 15-render-quadlets.sh does
# NOT process, and systemd does not expand env vars in User=/WorkingDirectory=/
# ConditionPathExists=/ExecStart=. So the values below are LITERALS that MIRROR
# mios.toml ([ports].opencode_gateway, [ai].opencode_*, [ai].agent_venv) — the
# same convention hermes-agent.service uses. If a value changes in mios.toml,
# update this unit to match.
#
# Skip cleanly if the opencode binary never landed (build-time fetch failure)
# rather than crash-loop; `mios update` re-runs 38-hermes-agent.sh (PHASE 2)
# to complete it.
ConditionPathExists=/usr/lib/mios/agents/opencode/bin/opencode

[Service]
Type=simple
# Agent-plane user consolidation: runs as the shared mios-ai user (matches
# hermes-agent.service + mios-agent-pipe.service), NOT root. mios-ai owns the
# /var/lib/mios state tree + reads /etc/mios config.
User=mios-ai
Group=mios-ai
# Gateway run cwd + HOME (state, owned by mios-ai; created by the
# mios-opencode-gateway.conf tmpfiles). Mirrors [ai].opencode_gateway_workdir.
WorkingDirectory=/var/lib/mios/opencode-gateway/work
Environment=HOME=/var/lib/mios/opencode-gateway/work
# SSOT mirrors: [ports].opencode_gateway, [ai].opencode_bin/opencode_model/
# opencode_provider/opencode_config, [ai].opencode_gateway_timeout_s.
Environment=MIOS_PORT_OPENCODE_GATEWAY=8633
Environment=MIOS_OPENCODE_HOST=127.0.0.1
Environment=MIOS_OPENCODE_BIN=/usr/lib/mios/agents/opencode/bin/opencode
Environment=MIOS_OPENCODE_MODEL=mios-opencode:latest
Environment=MIOS_OPENCODE_PROVIDER=local
Environment=MIOS_OPENCODE_CONFIG=/etc/mios/opencode/opencode.json
Environment=MIOS_OPENCODE_TIMEOUT_S=90
# Reuse the SHARED agent-plane venv interpreter (the shim is stdlib-only, but
# this matches the agent-pipe pattern + keeps ONE interpreter). This venv is the
# de-facto shared interpreter for hermes-agent, agent-pipe, and this gateway
# (mios.toml [ai].agent_venv = /usr/lib/mios/agents/.venv).
ExecStart=/usr/lib/mios/agents/.venv/bin/python3 -u /usr/lib/mios/agents/opencode-gateway/server.py
Restart=on-failure
RestartSec=10s
TimeoutStartSec=30s
# --- hardening (mirrors hermes-agent.service / mios-agent-pipe.service) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/mios /var/lib/mios
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
