# usr/lib/systemd/system/mios-freeipa-enroll.service
# AI-hint: Systemd unit that executes mios-freeipa-enroll.sh to perform zero-touch FreeIPA enrollment for WSL and non-containerized environments, triggered only if /etc/mios/ipa-enroll.env exists and /etc/ipa/default.conf is missing.
# AI-related: /etc/mios/ipa-enroll.env, /usr/libexec/mios/mios-freeipa-enroll.sh, mios-freeipa-enroll, network-online.target, multi-user.target
[Unit]
Description='MiOS' FreeIPA Zero-Touch Enrollment
ConditionPathExists=/etc/mios/ipa-enroll.env
ConditionPathExists=!/etc/ipa/default.conf
# Network-driven enrollment -- WSL deployments enrolling against a
# corporate IPA realm need this too. Gated by /etc/mios/ipa-enroll.env
# existence, so it stays a no-op unless explicitly configured.
ConditionVirtualization=|!container
ConditionVirtualization=|wsl
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mios/mios-freeipa-enroll.sh
# Hardening to protect the credential reading process
UMask=0077

[Install]
WantedBy=multi-user.target
