# usr/lib/systemd/system/mios-aios-refresh.timer
[Unit]
Description=Periodic MiOS AIOS refresh (SSOT role SYSTEMs + A2A fleet discovery)

[Timer]
# 2 min after boot (agent-pipe up by then), then every 15 min so the role
# SYSTEMs track catalog/SSOT changes and the peer list tracks the live fleet.
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
