# usr/lib/systemd/system/mios-forge-firstboot.service
[Unit]
Description='MiOS' Forge first-boot admin-bootstrap (Forgejo)
Documentation=https://github.com/mios-dev/MiOS/blob/main/etc/containers/systemd/mios-forge.container
After=mios-forge.service network-online.target
Wants=mios-forge.service network-online.target
ConditionPathExists=/etc/mios/install.env
ConditionPathExists=!/var/lib/mios/forge/.firstboot-done
# Forgejo bootstrap MUST run on WSL too (the localhost:3000 git origin
# is the heart of the MiOS self-replication loop on every shape). A
# plain `!container` would skip WSL because systemd reports virtualization
# as `wsl` -- a container variant. Use OR: "not a container OR is WSL".
ConditionVirtualization=|!container
ConditionVirtualization=|wsl

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/mios/install.env
ExecStart=/usr/libexec/mios/forge-firstboot.sh
TimeoutStartSec=600s

# Hardening: this service writes to a small set of paths plus calls
# 'podman exec' against the running mios-forge container. RestrictNamespaces
# and RestrictAddressFamilies were tried but break Podman's CRIU/conmon
# attach path on rootful container exec; we drop them and lean on the
# read-write path scoping + ProtectHome instead, which is sufficient for
# this script's actual surface area.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/mios/forge /var/lib/mios/forge /var/log/mios/forge /var/lib/containers /run/containers
PrivateTmp=yes
LockPersonality=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
