# usr/lib/systemd/system/mios-verify.service
# AI-hint: Systemd unit that executes /usr/libexec/mios-verify to perform a fs-verity cryptographic integrity audit of the root filesystem during the initrd boot sequence before mounting the sysroot.
# AI-related: mios-verify, ostree-prepare-root.service, sysroot.mount, initrd-root-fs.target
[Unit]
ConditionVirtualization=!container
Description='MiOS' Cryptographic Integrity Audit (fs-verity)
DefaultDependencies=no
After=ostree-prepare-root.service
Before=sysroot.mount
Before=ostree-prepare-root.service
ConditionPathExists=/sysroot/ostree/deploy

[Service]
Type=oneshot
# Execute the verification script to validate fs-verity hashes
ExecStart=/usr/libexec/mios-verify
RemainAfterExit=yes

[Install]
WantedBy=initrd-root-fs.target
