# usr/lib/systemd/system/mios-boot-diag.service
# AI-hint: Systemd unit file defining a one-shot service that executes /usr/libexec/mios-boot-diag to output boot status to the console and journal during the transition from basic.target to graphical.target.
# AI-related: mios-boot-diag, basic.target, graphical.target
[Unit]
Description='MiOS' Boot Diagnostic -- prints service status to console
# Read-only journal/systemctl introspection -- safe on every shape.
ConditionVirtualization=|!container
ConditionVirtualization=|wsl
After=basic.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/mios-boot-diag
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=graphical.target