# usr/lib/systemd/system/mios-ollama-firstboot.service
[Unit]
Description='MiOS' Ollama first-boot model pull
Documentation=https://ollama.com/
# Pull the default model set into /var/lib/ollama on first boot. We
# deliberately don't bake models into the image (build-time /var
# cleanup would wipe them). Sentinel /var/lib/mios/.ollama-firstboot-
# done marks completion; deleting it forces a re-pull on next start.
After=ollama.service network-online.target
Wants=ollama.service network-online.target
Requires=ollama.service
ConditionPathExists=!/var/lib/mios/.ollama-firstboot-done
# Run on bare-metal/VM AND on WSL. systemd treats `wsl` as a container
# variant, so a plain `!container` would skip every WSL deployment --
# but Ollama is a core MiOS feature on WSL too. The two-line ORed form
# below evaluates to "not a container OR is WSL".
ConditionVirtualization=|!container
ConditionVirtualization=|wsl

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mios/ollama-firstboot.sh
TimeoutStartSec=1800
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target default.target
