# usr/lib/systemd/system/mios-k3s-init.service
[Unit]
Description='MiOS' K3s Manifest Seeder
Before=k3s.service
ConditionPathExists=/usr/share/mios/k3s-manifests
# Pure file-copy oneshot -- no kernel/hardware coupling. Run on every
# shape including WSL (k3s.container itself stays !wsl-gated; this
# seeder is harmless if k3s ultimately doesn't start).
ConditionVirtualization=|!container
ConditionVirtualization=|wsl

[Service]
Type=oneshot
RemainAfterExit=yes

# Copy the immutable manifests into the persistent datastore
ExecStart=/bin/bash -c "cp -a /usr/share/mios/k3s-manifests/*.yaml /var/lib/rancher/k3s/server/manifests/ 2>/dev/null || true"
ExecStartPost=/bin/chmod 0600 /var/lib/rancher/k3s/server/manifests/*.yaml

[Install]
WantedBy=multi-user.target
