# usr/lib/systemd/system/mios-policy-arbiter.service
# AI-hint: Systemd unit for the WS-9 out-of-process HITL policy arbiter -- runs /usr/libexec/mios/mios-policy-arbiter (a stdlib loopback HTTP service) as the mios-ai user, answering the agent-pipe's HITL arbiter client with allow/deny verdicts decided by mios_arbiter over the operator policy (MIOS_ARBITER_* from install.env). Idle/no-op until [ai].hitl_arbiter_url points at it; default policy is allow-all so enabling it changes nothing until a deny-list/block-tier is set.
# AI-related: /usr/libexec/mios/mios-policy-arbiter, /usr/lib/mios/agent-pipe/mios_arbiter.py, /usr/lib/mios/agent-pipe/server.py, /usr/share/mios/mios.toml, mios-agent-pipe.service
# /usr/lib/systemd/system/mios-policy-arbiter.service
# 'MiOS' out-of-process HITL policy arbiter (WS-9). A second, operator-ownable
# opinion ON TOP of the in-process #62 HITL gate + WS-A9 PDP: the agent-pipe POSTs
# each high-risk (tier >= [ai].hitl_threshold) action here for an allow/deny
# verdict. Runs as mios-ai (least privilege); binds 127.0.0.1 only. The agent-pipe
# consults it ONLY when [ai].hitl_arbiter_url is set -> idle otherwise.
[Unit]
Description='MiOS' out-of-process HITL policy arbiter (WS-9)
After=network-online.target mios-agent-pipe.service
Wants=network-online.target
ConditionPathExists=/usr/libexec/mios/mios-policy-arbiter

[Service]
Type=simple
User=mios-ai
Group=mios-ai
Environment=MIOS_AGENT_PIPE_DIR=/usr/lib/mios/agent-pipe
EnvironmentFile=-/etc/mios/install.env
ExecStart=/usr/bin/python3 /usr/libexec/mios/mios-policy-arbiter
Restart=on-failure
RestartSec=15s
# Hardening: a stateless loopback decision service needs nothing but Python.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
