# usr/lib/systemd/system/mios-boot-diag.service
[Unit]
Description='MiOS' Boot Diagnostic -- prints service status to console
# Read-only journal/systemctl introspection -- safe on every shape.
ConditionVirtualization=|!container
ConditionVirtualization=|wsl
After=basic.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/mios-boot-diag
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=graphical.target