# usr/lib/systemd/system/mios-aios-refresh.timer
# AI-hint: Systemd timer that triggers mios-aios-refresh.service every 15 minutes to synchronize the Single Source of Truth (SSOT) and update the A2A fleet discovery peer list.
# AI-related: mios-aios-refresh, mios-aios-refresh.service, timers.target
[Unit]
Description=Periodic MiOS AIOS refresh (SSOT role SYSTEMs + A2A fleet discovery)

[Timer]
# 2 min after boot (agent-pipe up by then), then every 15 min so the role
# SYSTEMs track catalog/SSOT changes and the peer list tracks the live fleet.
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
