# usr/lib/systemd/system/mios-wsl-firstboot.service
# AI-hint: Systemd unit that executes /usr/libexec/mios/wsl-firstboot to perform one-time initialization tasks specifically for WSL2 instances if the first-boot flag is not yet set.
# AI-related: /usr/libexec/mios/wsl-firstboot, systemd-user-sessions.service, local-fs.target, multi-user.target
[Unit]
Description='MiOS' WSL2 First Boot Initialization
Documentation=https://github.com/MiOS-DEV/MiOS
ConditionVirtualization=wsl
ConditionPathExists=!/var/lib/mios/.wsl-firstboot-done
After=local-fs.target
Before=systemd-user-sessions.service multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mios/wsl-firstboot

[Install]
WantedBy=multi-user.target
