# usr/lib/systemd/system/mios-freeipa-enroll.service
[Unit]
Description='MiOS' FreeIPA Zero-Touch Enrollment
ConditionPathExists=/etc/mios/ipa-enroll.env
ConditionPathExists=!/etc/ipa/default.conf
# Network-driven enrollment -- WSL deployments enrolling against a
# corporate IPA realm need this too. Gated by /etc/mios/ipa-enroll.env
# existence, so it stays a no-op unless explicitly configured.
ConditionVirtualization=|!container
ConditionVirtualization=|wsl
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mios/mios-freeipa-enroll.sh
# Hardening to protect the credential reading process
UMask=0077

[Install]
WantedBy=multi-user.target
