# usr/lib/systemd/system/mios-nvidia-cdi.service
# AI-hint: Systemd unit that executes nvidia-ctk to generate the NVIDIA CDI specification file at /etc/cdi/nvidia.yaml, enabling GPU passthrough for rootless Podman containers.
# AI-related: nvidia-persistenced.service, multi-user.target
[Unit]
Description=Generate NVIDIA CDI spec for rootless Podman GPU
After=nvidia-persistenced.service
ConditionPathExists=/usr/bin/nvidia-ctk
ConditionVirtualization=!vm
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nvidia-ctk cdi generate --output=/etc/cdi/nvidia.yaml
ExecStartPost=/usr/bin/nvidia-ctk cdi list
[Install]
WantedBy=multi-user.target