# usr/lib/systemd/system/mios-delegation-prefilter.service
[Unit]
# MiOS delegation prefilter -- HTTP forwarder that sits between Open
# WebUI (and any other OpenAI-compat client) and Hermes-Agent's
# /v1/chat/completions, force-routing fan-outable user prompts through
# delegate_task by injecting tool_choice on the request.
#
# Why a sidecar instead of a Hermes plugin: Hermes's plugin discovery
# is documented but the hook signatures are async/sync mixed, and an
# async/sync wrapper mismatch on the wrong handler crashes EVERY tool
# call ('object of type coroutine has no len()' -- operator-confirmed
# regression 2026-05-15). A separate-process forwarder is:
#   * deterministic (heuristic-driven, not model-driven)
#   * isolated from Hermes upgrades
#   * trivially disablable (stop the unit, point OWUI back at :8642)
# Cost: one extra hop per request (~1ms localhost overhead).
Description='MiOS' delegation prefilter -- forces delegate_task on fan-out prompts
Documentation=file:///usr/libexec/mios/mios-delegation-prefilter
After=hermes-agent.service network.target
Wants=hermes-agent.service
ConditionPathExists=/usr/lib/mios/hermes-agent/.venv/bin/python3

[Service]
Type=exec
ExecStart=/usr/lib/mios/hermes-agent/.venv/bin/python3 /usr/libexec/mios/mios-delegation-prefilter
Restart=on-failure
RestartSec=3s
User=mios-ai
Group=mios-ai
Environment=PYTHONUNBUFFERED=1
# Operator architecture 2026-05-16: MiOS-Agent (the OWUI face) is a
# small CPU model that REFINES the user prompt with full MiOS context
# BEFORE the heavier hermes orchestrator sees it. The refiner is
# mios-sys-agent (~2.7 GB, runs on CPU). Enabling MIOS_PREFILTER_REFINE
# here means every prompt OWUI sends actually gets the small-CPU-model
# rewrite before hermes -- not a dumb proxy. The refined prompt + the
# original both end up at hermes; the refiner output is also written
# to /var/lib/mios/delegation-prefilter/latest.json so the OWUI sidecar
# filter can surface it as a "🧠 refined: ..." status event.
# Refinement is the OWUI MiOS-Agent PIPE's job (it owns the CPU
# refine step inline so OWUI sees the progress emits). The prefilter
# only needs to force delegate_task + auth-passthrough now. Keep the
# env-var schema for back-compat; set =0 by default. Operators with
# non-OWUI clients (curl scripts, MCP) that need server-side refine
# can flip this to =1 via a drop-in.
Environment=MIOS_PREFILTER_REFINE=0
Environment=MIOS_SYS_AGENT_MODEL=mios-sys-agent
Environment=MIOS_SYS_AGENT_BASE=http://127.0.0.1:11434/v1
Environment=MIOS_SYS_AGENT_ENDPOINT=http://127.0.0.1:11434
# Cold-load of the ~4 GB CPU refiner takes 60-180s on WSL2 disks
# (3+ min observed on the operator's slow-disk testbed); on bare-
# metal NVMe it's 10-30s. After cold load, mios-sys-agent generates
# ~50 tok/s on a 12-core CPU, so a 600-token refine completes in
# ~12s warm. 300s total covers cold + warm + the elaborate
# modelfile SYSTEM prompt (~1500 tokens of reasoning template).
# Combined with the OLLAMA_KEEP_ALIVE=-1 drop-in the cold load
# happens ONCE per boot and never again.
Environment=MIOS_PREFILTER_REFINE_TIMEOUT=300
# Heuristic + upstream are operator-tunable via /etc/mios/secrets.env
# or a drop-in:
#   Environment=MIOS_PREFILTER_LISTEN_PORT=8641
#   Environment=MIOS_PREFILTER_UPSTREAM=http://127.0.0.1:8642

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
# Refiner writes the sidecar state file here so the OWUI mios_sidecar
# Filter can emit "🧠 mios-sys-agent refined prompt" status events.
# Without this RW exception ProtectSystem=strict blocks the write
# silently and the operator never sees the refinement happen.
ReadWritePaths=/var/lib/mios/delegation-prefilter
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
