# usr/lib/systemd/system/mios-ha-node.target
# AI-hint: Defines the mios-ha-node.target unit to identify and configure a node as a High Availability cluster member, ensuring it conflicts with desktop/headless modes and requires corosync/pacemaker services.
# AI-related: mios-ha-node, mios-desktop, mios-headless, mios-k3s-master, corosync.service, pacemaker.service, multi-user.target, mios-desktop.target, mios-headless.target, mios-k3s-master.target
[Unit]
Description='MiOS' HA Cluster Node Role
Requires=multi-user.target corosync.service pacemaker.service
After=multi-user.target
Conflicts=mios-desktop.target mios-headless.target mios-k3s-master.target
AllowIsolate=yes
[Install]
WantedBy=multi-user.target
