# usr/lib/systemd/system/mios-skills-miner.timer
# /usr/lib/systemd/system/mios-skills-miner.timer
# Phase C.2 of the AgentOS roadmap: cadence for the background
# skill miner. Interval lifted to mios.toml [skills].
# mine_interval_minutes (default 60). Operator override:
#   sudo systemctl edit mios-skills-miner.timer
#   [Timer]
#   OnUnitActiveSec=30min
#
# Disabled by default; operator opts in (or it inherits enablement
# from the configurator HTML "Skills mining" toggle which maps to
# [skills].enable). The .service ConditionPathExists guard means a
# stripped-down deployment with the libexec script absent skips
# silently.

[Unit]
Description='MiOS' Phase C.2 skill-miner cadence (sequential pattern mining)
Documentation=file:///usr/libexec/mios/mios-skills

[Timer]
OnBootSec=10min
OnUnitActiveSec=60min
AccuracySec=2min
Unit=mios-skills-miner.service
Persistent=true

[Install]
WantedBy=timers.target
