# usr/lib/systemd/system/cockpit.service.d/10-mios-container.conf
# /usr/lib/systemd/system/cockpit.service.d/10-mios-container.conf
#
# Relax cockpit's namespace-isolation hardening so the unit comes up inside
# containers and WSL2. Cockpit's upstream unit (in cockpit-ws) configures
# cockpit-certificate-ensure with PrivateNetwork=yes, PrivateIPC=yes,
# PrivateMounts=yes, ProtectHostname=yes, ProtectControlGroups=yes,
# ProtectKernelTunables=yes -- none of which the kernel allows when systemd
# is running without CAP_SYS_ADMIN over its parent namespaces, which is the
# normal case when MiOS runs as a podman-machine, distrobox, or nested
# container.
#
# Symptom on the 2026-05-05 dev VM run:
#
#     cockpit.service: PrivateNetwork=yes is configured, but the kernel
#         does not support or we lack privileges for network namespace,
#         proceeding without.
#     cockpit.service: PrivateIPC=yes is configured, but IPC namespace
#         setup failed, ignoring: Operation not permitted
#     cockpit.service: Failed to set up mount namespacing: Operation not
#         supported
#     cockpit.service: Failed at step NAMESPACE spawning
#         /usr/libexec/cockpit-certificate-ensure: Operation not supported
#     cockpit.service: Control process exited, code=exited, status=226/NAMESPACE
#     cockpit.service: Failed with result 'exit-code'.
#
# The retries hit cockpit.socket's TriggerLimitBurst (5 in 10 s) and Cockpit
# falls over completely. Bluefin and Aurora ship the equivalent drop-in
# unconditionally; cockpit-certificate-ensure is a tiny binary whose only
# job is to touch a TLS cert under /etc/cockpit/ws-certs.d/, so the
# hardening it ships with is essentially symbolic and the relaxation cost is
# negligible. The actual long-lived cockpit.service runtime is unaffected
# because cockpit's privsep happens via cockpit-bridge -> cockpit-session
# inside cockpit-ws's own user namespace.

[Service]
PrivateNetwork=no
PrivateIPC=no
PrivateMounts=no
ProtectHostname=no
ProtectControlGroups=no
ProtectKernelTunables=no
