# usr/lib/systemd/system/cockpit.service.d/10-mios-container.conf
# /usr/lib/systemd/system/cockpit.service.d/10-mios-container.conf
#
# Relax cockpit's namespace-isolation hardening so the unit comes up inside
# containers and WSL2. Cockpit's upstream unit (in cockpit-ws) configures
# cockpit-certificate-ensure with the full systemd hardening surface --
# PrivateNetwork, PrivateIPC, PrivateMounts, PrivateTmp, ProtectSystem,
# ProtectHome, ProtectKernelTunables, ProtectKernelModules, ProtectClock,
# RestrictNamespaces, NoNewPrivileges, plus a SystemCallFilter -- none of
# which the kernel allows when systemd is running without CAP_SYS_ADMIN
# over its parent namespaces. That is the normal case when MiOS runs as
# a podman-machine, distrobox, WSL2 distro, or nested container.
#
# Symptom on the 2026-05-05 / 2026-05-06 WSL2 boot:
#
#     cockpit.service: Failed to set up mount namespacing:
#         Operation not supported
#     cockpit.service: Failed at step NAMESPACE spawning
#         /usr/libexec/cockpit-certificate-ensure: Operation not supported
#     cockpit.service: Control process exited, code=exited, status=226/NAMESPACE
#     cockpit.service: Failed with result 'exit-code'.
#     cockpit.socket: Failed with result 'service-start-limit-hit'.
#
# Even with PrivateMounts=no set, ANY of {PrivateTmp, ProtectSystem,
# ProtectHome, ProtectKernelTunables, ProtectControlGroups,
# RestrictNamespaces, ProtectKernelModules, ProtectClock, ProtectKernelLogs,
# ReadWritePaths, ReadOnlyPaths, InaccessiblePaths} causes systemd to clone
# a mount namespace at exec time, and that clone fails with EOPNOTSUPP on
# any host where the parent namespace is unprivileged. Below we explicitly
# neutralize ALL of them so cockpit-certificate-ensure can exec without any
# CLONE_NEWNS attempt at all.
#
# The relaxation cost is negligible. cockpit-certificate-ensure is a tiny
# binary whose only job is to ensure a TLS cert exists under
# /etc/cockpit/ws-certs.d/ -- it doesn't touch the kernel surface.
# Cockpit's actual long-lived runtime privsep happens later, via
# cockpit-bridge -> cockpit-session inside cockpit-ws's own user namespace,
# which is unaffected by this drop-in. Bluefin / Aurora / Universal Blue
# ship the equivalent drop-in unconditionally for the same reason.

[Service]
# Namespace isolation -- every directive that triggers CLONE_NEWNS:
PrivateNetwork=no
PrivateIPC=no
PrivateMounts=no
PrivateTmp=no
PrivateDevices=no
PrivateUsers=no
# ProtectSystem= and ProtectHome= remount /usr / /etc / /home read-only via
# a private mount namespace. Both must be off for the EOPNOTSUPP path.
ProtectSystem=no
ProtectHome=no
ProtectHostname=no
ProtectClock=no
ProtectKernelTunables=no
ProtectKernelModules=no
ProtectKernelLogs=no
ProtectControlGroups=no
ProtectProc=default
ProcSubset=all
# Process / capability hardening -- not strictly mount-namespace, but the
# certificate-ensure helper doesn't need them and they confuse some
# minimal containers (LXC, podman with --privileged=false).
NoNewPrivileges=no
RestrictNamespaces=no
RestrictRealtime=no
RestrictSUIDSGID=no
LockPersonality=no
MemoryDenyWriteExecute=no
# Reset list-typed directives so any inherited ReadWritePaths=/...
# doesn't smuggle a private mount back in.
SystemCallFilter=
SystemCallArchitectures=
RestrictAddressFamilies=
DeviceAllow=
ReadWritePaths=
ReadOnlyPaths=
InaccessiblePaths=
