# usr/lib/systemd/system/cockpit.service.d/10-mios-wsl2.conf
[Service]
# WSL2: give cockpit extra startup time (slower I/O path vs bare metal).
TimeoutStartSec=90
