# usr/lib/systemd/system/cockpit.service.d/10-mios-wsl2.conf
[Unit]
# cockpit-ws uses systemd-logind varlink for session management.
# In WSL2, logind retries at boot -- wait for it to be fully active.
After=systemd-logind.service dbus-daemon-wsl.service

[Service]
# Give cockpit-ws extra time to initialize TLS certs on first run
TimeoutStartSec=90
