# usr/lib/systemd/system/cockpit.service.d/10-mios-wsl2.conf
[Unit]
ConditionVirtualization=!wsl

[Service]
TimeoutStartSec=90
