# usr/lib/systemd/system/cockpit.service.d/10-mios-wsl2.conf
# AI-hint: Configures the cockpit.service unit to allow a 90s TimeoutStartSec specifically for WSL2 environments to accommodate slower I/O paths during startup.
# AI-related: cockpit.service
[Service]
# WSL2: give cockpit extra startup time (slower I/O path vs bare metal).
TimeoutStartSec=90
