# usr/lib/systemd/system/systemd-resolved-varlink.socket.d/10-mios-skip-container.conf
# /usr/lib/systemd/system/systemd-resolved-varlink.socket.d/10-mios-skip-container.conf
#
# Skip the resolved varlink socket inside containers where resolved is
# gated off. Same failure pattern as 10-mios-skip-container.conf for the
# logind varlink socket: resolved.service has
# ConditionVirtualization=!container, so its varlink socket activates
# repeatedly to a service that gets skipped, hits TriggerLimitBurst, and
# leaves a permanently-failed socket in the unit list. Adding the same
# Condition* directives to the socket prevents the activation churn.

[Unit]
ConditionVirtualization=!container
