# usr/lib/systemd/system/upower.service.d/10-mios-wsl2.conf
# AI-hint: Disables the upower service in WSL2 environments by checking for the wsl_vm_type to prevent SIGABRT crashes caused by missing hardware power_supply interfaces.
[Unit]
# upower probes /sys/class/power_supply and DBus hardware interfaces that
# don't exist in WSL2, causing an immediate SIGABRT. Skip it entirely.
ConditionVirtualization=!wsl
