# usr/lib/systemd/system/upower.service.d/10-mios-wsl2.conf
[Unit]
# upower probes /sys/class/power_supply and DBus hardware interfaces that
# don't exist in WSL2, causing an immediate SIGABRT. Skip it entirely.
ConditionVirtualization=!wsl
