# usr/lib/systemd/system/cockpit-wsinstance-http.socket.d/10-mios-container.conf
# /usr/lib/systemd/system/cockpit-wsinstance-http.socket.d/10-mios-container.conf
#
# Pair drop-in to cockpit-wsinstance-socket-user.service.d/10-mios-container.conf.
# After the socket-user service drops its namespace requirements, the
# child socket units (cockpit-wsinstance-http.socket and
# cockpit-wsinstance-https-factory.socket) inherit DynamicUser /
# PrivateMounts from cockpit's stock unit and fail their own
# `Failed to listen on cockpit-wsinstance-http.socket` step inside a
# container or podman-machine WSL distro. Strip the same namespace
# directives so the .sock binds succeed.
#
# KEEP the stock SocketUser=root + SocketGroup=cockpit-wsinstance-socket
# (do NOT override to cockpit-ws:cockpit-ws). Operator-flagged 2026-05-10:
# overriding to cockpit-ws:cockpit-ws caused
#     cockpit-tls: connect(http.sock) failed: Permission denied
# because cockpit.service runs as User=cockpit-systemd-service /
# Group=cockpit-wsinstance-socket, which has no group-rw access to a
# cockpit-ws-owned socket.

[Unit]
ConditionVirtualization=container

[Socket]
DynamicUser=no
PrivateTmp=no
PrivateMounts=no
PrivateDevices=no
PrivateUsers=no
ProtectHome=no
ProtectSystem=no
RestrictNamespaces=no
