# usr/lib/systemd/system/mios-delegation-prefilter.service.d/10-sys-agent-refine.conf
# /usr/lib/systemd/system/mios-delegation-prefilter.service.d/
#   10-sys-agent-refine.conf
#
# Enable the mios-sys-agent refinement stage in the delegation
# prefilter. Operator directive 2026-05-16: "mios-sys-agent should
# dispatch to the current MiOS-Hermes with the tasks to delegate --
# mios-sys-agent simply refines the prompt with the provided template
# key points to other MiOS-Agent(s)".
#
# When enabled: every chat-completions request to :8641 first goes
# to mios-sys-agent (the Ollama tag derived from qwen3:0.6b-cpu with
# the reasoning template + MiOS environment map baked into SYSTEM).
# Sys-agent emits a structured handoff (USER REQUEST / INTENT /
# CONSTRAINTS / MIOS CONTEXT / PLAN). The handoff replaces the
# user message; hermes-agent processes the enriched prompt.
#
# Disable by removing this drop-in or setting MIOS_PREFILTER_REFINE=0.

[Service]
# Refinement RE-ENABLED 2026-05-16 -- mios-sys-agent now backed by
# qwen3.5:2b on GPU (~4 GB VRAM, ~200-500 ms per refinement when
# warm). Prior CPU-pinned 30-60s latency is gone with GPU residency.
# firstboot pre-warm keeps it loaded with keep_alive=-1.
#
# Architecture per operator directive 2026-05-16:
#   operator types in OWUI
#     -> prefilter receives chat completion
#     -> mios-sys-agent (qwen3.5:2b GPU) applies reasoning template
#        + restructures prompt into INTENT / CONSTRAINTS / MIOS
#        CONTEXT / PLAN handoff
#     -> hermes (qwen3-coder:30b GPU) executes against refined prompt
#     -> response streams back to operator
Environment=MIOS_PREFILTER_REFINE=1
Environment=MIOS_SYS_AGENT_MODEL=mios-sys-agent
Environment=MIOS_SYS_AGENT_ENDPOINT=http://127.0.0.1:11434
# Model-name rewrite -- OWUI shows "MiOS-Agent" in dropdown; the
# prefilter translates to "hermes-agent" on the way in.
Environment=MIOS_PREFILTER_DISPLAY_MODEL=MiOS-Agent
# 45 s gives ample headroom for cold-load (rare with pre-warm); warm
# calls complete in <1s.
Environment=MIOS_PREFILTER_REFINE_TIMEOUT=45
