# usr/lib/systemd/system/mios-delegation-prefilter.service.d/10-sys-agent-refine.conf
# /usr/lib/systemd/system/mios-delegation-prefilter.service.d/
#   10-sys-agent-refine.conf
#
# DEPRECATED 2026-05-17. Refinement moved from the prefilter to the
# OWUI MiOS-Agent pipe per operator architecture restatement:
# "OPEN WEBUI'S MIOS_AGENT OPERATES ON THE CPU MODEL (PROCESSED ON
# THE CPU) -- IS QUICKLY REFINING THE USERS PROMPTS". The pipe at
# /usr/share/mios/owui/pipes/mios_agent_pipe.py now owns the refine
# step inline; the prefilter is back to its original purpose of
# forcing delegate_task on every user turn + auth passthrough.
#
# This drop-in is preserved (rather than deleted) so it survives
# `dnf reinstall` / `bootc upgrade` migration and explicitly
# REINFORCES prefilter-side refine OFF -- and timeout HIGH for the
# corner case where the operator re-enables server-side refine.
#
# Operators wanting server-side refinement for non-OWUI clients
# (raw curl, MCP, etc.) should:
#   1. Set Environment=MIOS_PREFILTER_REFINE=1 here
#   2. Set REFINE_ENABLED=False in the pipe's valves (OWUI admin UI)
#   ...so refinement happens in EXACTLY ONE place (not both).

[Service]
Environment=MIOS_PREFILTER_REFINE=0
Environment=MIOS_SYS_AGENT_MODEL=mios-sys-agent
Environment=MIOS_SYS_AGENT_ENDPOINT=http://127.0.0.1:11434
Environment=MIOS_PREFILTER_DISPLAY_MODEL=MiOS-Agent
Environment=MIOS_PREFILTER_REFINE_TIMEOUT=300
