# usr/lib/systemd/system/mios-delegation-prefilter.service.d/10-sys-agent-refine.conf
# /usr/lib/systemd/system/mios-delegation-prefilter.service.d/
#   10-sys-agent-refine.conf
#
# Enable the mios-sys-agent refinement stage in the delegation
# prefilter. Operator directive 2026-05-16: "mios-sys-agent should
# dispatch to the current MiOS-Hermes with the tasks to delegate --
# mios-sys-agent simply refines the prompt with the provided template
# key points to other MiOS-Agent(s)".
#
# When enabled: every chat-completions request to :8641 first goes
# to mios-sys-agent (the Ollama tag derived from qwen3:0.6b-cpu with
# the reasoning template + MiOS environment map baked into SYSTEM).
# Sys-agent emits a structured handoff (USER REQUEST / INTENT /
# CONSTRAINTS / MIOS CONTEXT / PLAN). The handoff replaces the
# user message; hermes-agent processes the enriched prompt.
#
# Disable by removing this drop-in or setting MIOS_PREFILTER_REFINE=0.

[Service]
# Refinement DISABLED 2026-05-16 (operator-observed: CPU-pinned
# qwen3:0.6b-cpu + reasoning-template SYSTEM prompt = 30-60s per
# refine call; unacceptable for the chat path). The right home for
# the sys-agent refinement role is an OWUI Pipe Function (runs in
# OWUI, can stream emitter events, can hit OWUI RAG/knowledge/
# memory APIs natively, can call hermes as the worker). Follow-up
# commit will build that pipe; until then, REFINE=0 keeps the chat
# path snappy.
Environment=MIOS_PREFILTER_REFINE=0
Environment=MIOS_SYS_AGENT_MODEL=mios-sys-agent
Environment=MIOS_SYS_AGENT_ENDPOINT=http://127.0.0.1:11434
# Model-name rewrite stays ON -- OWUI shows "mios-sys-agent" in the
# dropdown; the prefilter translates to "hermes-agent" on the way
# in. From the operator's perspective they're chatting with mios-
# sys-agent (the name); the actual worker is hermes for now.
Environment=MIOS_PREFILTER_DISPLAY_MODEL=MiOS-Agent
Environment=MIOS_PREFILTER_REFINE_TIMEOUT=30
