# usr/lib/systemd/system/systemd-resolved.service.d/10-mios-container-gate.conf
# AI-hint: Configures systemd-resolved to disable container-specific logic when running on bare metal, ensuring consistent DNS resolution for MiOS services outside of OCI containers.
[Unit]
# 'MiOS': Redundant in OCI containers (handled by Podman/Docker)
ConditionVirtualization=!container
