# usr/lib/systemd/system/systemd-resolved.service.d/10-mios-container-gate.conf
[Unit]
# MiOS: Redundant in OCI containers (handled by Podman/Docker)
ConditionVirtualization=!container
