# usr/lib/systemd/system/coreos-ignition-firstboot-complete.service.d/10-mios-wsl2.conf
# AI-hint: Overrides the coreos-ignition-firstboot-complete.service to skip Ignition provisioning checks when running in a WSL2 environment where the boot partition label is unavailable.
# AI-related: coreos-ignition-firstboot-complete.service
[Unit]
# No Ignition provisioning in WSL2; /dev/disk/by-label/boot doesn't exist
ConditionVirtualization=!wsl
