# usr/lib/systemd/system/coreos-ignition-firstboot-complete.service.d/10-mios-wsl2.conf
[Unit]
# No Ignition provisioning in WSL2; /dev/disk/by-label/boot doesn't exist
ConditionVirtualization=!wsl
