# usr/lib/systemd/system/greenboot-healthcheck.service.d/10-network-online.conf
# AI-hint: Configures the greenboot-healthcheck service to delay execution until the network-online.target is reached, ensuring the health check agent has active connectivity to remote endpoints.
# AI-related: network-online.target
[Unit]
After=network-online.target
Wants=network-online.target
