# usr/lib/systemd/system/systemd-logind-varlink.socket.d/10-mios-skip-container.conf
# /usr/lib/systemd/system/systemd-logind-varlink.socket.d/10-mios-skip-container.conf
#
# Skip the logind varlink socket inside containers/WSL where logind itself
# is gated off. Without this drop-in the socket activates, the activated
# logind.service is skipped due to its own ConditionVirtualization=!wsl,
# the activation fails, the socket retries until TriggerLimitBurst is hit:
#
#     systemd[1]: systemd-logind-varlink.socket: Trigger limit hit,
#         refusing further activation.
#     systemd[1]: systemd-logind-varlink.socket: Failed with result
#         'trigger-limit-hit'.
#
# Failed-units count then includes the socket on every boot. Adding
# ConditionVirtualization to the [Unit] section makes systemd skip the
# socket entirely on container/WSL hosts, so the failure cascade never
# starts. Upstream applies the same condition to logind.service itself;
# this just propagates it to the matching varlink socket.

[Unit]
ConditionVirtualization=!container
ConditionVirtualization=!wsl
