# etc/systemd/user/localsearch-control-3.service.d/10-mios-allow-system-uid.conf
# /etc/systemd/user/localsearch-3.service.d/10-mios-allow-system-uid.conf
#
# MiOS dev VM creates the `mios` user with UID 992 (system class
# range, 0-999). Upstream localsearch-3 has ConditionUser=!@system
# so it refuses to start for system users. Without this drop-in,
# nautilus + gnome-files dbus-activate Tracker3.Miner.Files,
# localsearch-3 fails the condition check, dbus returns
# ServiceUnknown, the GTK app null-derefs in its bookmark loader.
#
# Future: ideally the mios user gets UID 1000+ at install time
# (mios-bootstrap useradd -u 1000) and this drop-in becomes
# unnecessary. Until then this is the live workaround.
[Unit]
ConditionUser=

