# usr/lib/systemd/system/nvidia-powerd.service.d/10-mios-virt-gate.conf
[Unit]
# MiOS: Skip in WSL2 and OCI containers for stability and performance
ConditionVirtualization=!wsl
ConditionVirtualization=!container
