# usr/share/containers/systemd/users/mios-coderun-sandbox@.container
# AI-hint: A Quadlet-managed systemd user unit defining a podman container for executing untrusted code in a dry-run sandbox, providing a secure, isolated environment for agents to test scripts before host execution.
# AI-related: /etc/mios/containers/coderun-seccomp.json, mios-coderun-sandbox, mios-coderun-session, mios-ai, mios-ai-firstboot, mios-ceph, mios-k3s, mios-coderun, default.target, network-online.target
[Unit]
# /usr/share/containers/systemd/users/mios-coderun-sandbox@.container
#
# CODE-RUN sandbox (Quadlet, templated). NOT the agent runtime --
# MiOS agents (Hermes, sys-agent, opencode, micro-LLMs) ALWAYS install
# directly to the host root per operator architecture. This container
# is the dry-run / code-test boundary that the agent (or operator)
# dispatches CODE into when they want to verify behavior before
# touching the host -- e.g. "run this generated build script in a
# disposable namespace and report exit code + stdout".
#
# USER-SCOPE Quadlet (operator 2026-06-12): this unit is started via
# `systemctl --user start mios-coderun-sandbox@<id>.service` by
# mios-coderun-session, and it uses %t (XDG_RUNTIME_DIR) + WantedBy=
# default.target -- i.e. it is a USER unit. It MUST live under a user
# Quadlet search path (here: /usr/share/containers/systemd/users/) so the
# podman USER generator emits it; when it was in the SYSTEM path
# (/usr/share|/etc/containers/systemd/) only a SYSTEM unit was generated
# and `systemctl --user` reported "unit not found" -> run_sandboxed_code
# failed with a sandbox/permission error. The mios-ai service account also
# needs linger + a subuid/subgid range (UserNS=keep-id) -- both provisioned
# by mios-ai-firstboot.
#
# Operator directive 2026-05-16: "MiOS-Agents live on the local MiOS
# systems -- Containers are for servers, applications, hosts, etc-etc
# (we will add proper sandboxing later (MiOS AI Agents ALWAYS install
# to the core system/host(s) root!!! Sandboxing can be later
# implemented with proper dry-running of code in these sandboxes or
# testing)".
#
# Start a dry-run sandbox with:
#
#   systemctl --user start mios-coderun-sandbox@<run-id>.service
#
# Dispatch code into it:
#
#   podman exec -i mios-coderun-sandbox-<run-id> \
#       /usr/local/bin/exec-init <cmd>
#
# The exec-init wrapper applies the per-process Landlock domain;
# the Quadlet provides namespace + seccomp + cgroup + network deny.
# Defense in depth: if any single layer fails, the others hold.

Description=MiOS code-run sandbox (%i) -- dry-run / test boundary, NOT agent runtime
After=network-online.target

[Container]
Image=localhost/mios-coderun-sandbox:latest
ContainerName=mios-coderun-sandbox-%i

# UNPRIVILEGED-QUADLETS (Law 6). UserNS=keep-id already pins the
# container "root" to whatever UID the unit runs as -- but the
# postcheck Law-6 validator REQUIRES an explicit User= directive on
# every Quadlet (no exceptions outside the documented allowlist:
# mios-ceph, mios-k3s). Set User=root because:
#   * the sandbox spawns container PID-1 as "root" inside the user
#     namespace (which keep-id maps to the host invoker, NOT real
#     uid 0); the inside-container root has NO host privileges
#   * the sandbox's whole defense-in-depth (Network=none, ReadOnly,
#     DropCapability=ALL, seccomp, Landlock) assumes container-root
#   * setting it to anything else breaks the keep-id mapping
# This satisfies the Law-6 grep while preserving the sandbox's
# rootless-on-the-host posture.
User=root

# UserNS=keep-id: container "root" is the host UID running this unit.
# Files written to bind-mounted /work appear as the operator's UID
# on the host -- no chown dance after the run.
UserNS=keep-id

# Network=none -- HARD network deny at the namespace level. Combined
# with the seccomp socket(AF_INET) ERRNO rule, the sandbox cannot
# reach the network even if it wanted to.
Network=none

# ReadOnly=true -- the container root filesystem is RO. Writes only
# allowed to the tmpfs mounts + the bind-mounted /work below.
ReadOnly=true

# VolatileTmp=true -- /tmp is its own tmpfs, wiped on container stop.
VolatileTmp=true

# Tmpfs mounts. /tmp for general scratch; /work-upper for overlayfs
# upperdir if the caller wants per-run COW on top of read-only bind.
# noexec on /work-upper just for paranoia.
Tmpfs=/tmp:rw,size=1G,mode=1777,nosuid,nodev
Tmpfs=/work-upper:rw,size=2G,mode=0755,nosuid,nodev,noexec

# No new privileges -- prctl(PR_SET_NO_NEW_PRIVS). Required for
# seccomp + Landlock from within an unprivileged container, and
# also a hard barrier against setuid escalation.
NoNewPrivileges=true

# Drop ALL capabilities. The code being run does NOT need
# CAP_NET_BIND_SERVICE (unprivileged ports), CAP_SYS_PTRACE
# (no debugger), CAP_SYS_ADMIN (no mount/umount), etc.
DropCapability=ALL

# SELinux container_t with per-instance MCS category -- only this
# container can read the :Z-labeled bind mount.
SecurityLabelType=container_t

# Hardened seccomp profile (allowlist; AF_UNIX-only sockets;
# clone3 -> ENOSYS; escalation syscalls -> ENOSYS).
SeccompProfile=/etc/mios/containers/coderun-seccomp.json

# Workspace bind mount, with per-instance SELinux label (:Z). The
# path is template-substituted: %i is the systemd instance specifier,
# used as both the run id AND the workspace directory name under
# /var/home/mios/coderuns/. (Operators wanting a different layout
# edit this Volume= line; future iteration reads from mios.toml
# [paths.coderun_workspace_root] per TOML-first invariant.)
Volume=/var/home/mios/coderuns/%i:/work:Z

# Unix socket for caller <-> sandbox RPC if running an in-container
# exec-server instead of `podman exec` from outside. %t = XDG_RUNTIME_DIR.
Volume=%t/mios-coderun-%i.sock:/run/coderun.sock

# Resource caps. Per the FS+OS Control guide §5 cgroups v2 section.
PodmanArgs=--cpus=4 \
           --memory=8g \
           --memory-swap=8g \
           --pids-limit=1024 \
           --ulimit nofile=4096 \
           --ulimit nproc=512 \
           --ulimit fsize=2147483648

[Service]
Restart=on-failure
RestartSec=2s
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=default.target
