# usr/lib/systemd/system/podman.socket.d/10-mios-operator-access.conf
# /usr/lib/systemd/system/podman.socket.d/10-mios-operator-access.conf
#
# Grant the OPERATOR's user access to the ROOTFUL podman API socket so SSH
# container managers (Serverbox, Termius, Portainer-likes) and the operator's
# own `podman` see + manage the real (rootful) MiOS containers. Default socket
# mode is 0600 root:root -> only root; a non-root manager app connecting as the
# operator is denied and shows nothing. Operator 2026-05-23: "Serverbox iphone
# app STILL doesn't show containers".
#
# Scope: group `mios` (the login operator, gid 1000) only -- NOT world. The
# hardened agent user (mios-ai) is deliberately NOT in this group; agents keep
# the READ-ONLY snapshot path (container_status verb) -- a network-facing
# service should not hold rootful container CONTROL. This grant is for the
# system OWNER's personal management surface, which inherently needs the API.
[Socket]
SocketMode=0660
SocketGroup=mios
