# usr/lib/systemd/system/crowdsec.service.d/10-bare-metal-only.conf
[Unit]
# MiOS: Skip this service in VMs/containers — bare metal only
ConditionVirtualization=no
