# usr/lib/systemd/system/systemd-machined.service.d/wsl2-optional.conf
# AI-hint: Configures systemd-machined to ignore cgroup-related failures in WSL2 environments by making the service non-fatal when the ConditionVirtualization=!wsl check fails.
[Unit]
# Make machined non-fatal in WSL2 (it needs cgroup features WSL2 lacks)
ConditionVirtualization=!wsl
