# usr/lib/systemd/system/systemd-machined.service.d/wsl2-optional.conf
[Unit]
# Make machined non-fatal in WSL2 (it needs cgroup features WSL2 lacks)
ConditionVirtualization=!wsl
