# usr/lib/systemd/system/mios-cdi-detect.service.d/50-gpu-passthrough.conf
# AI-hint: Configures the mios-cdi-detect service to execute mios-gpu-passthrough, automatically reconciling AddDevice drop-ins for background-AI Quadlets to ensure consistent iGPU/GPU visibility across all AI agents.
# AI-related: /usr/libexec/mios/mios-gpu-passthrough, mios-cdi-detect, mios-gpu-passthrough, mios-open-webui, mios-cdi-detect.service
# /usr/lib/systemd/system/mios-cdi-detect.service.d/50-gpu-passthrough.conf
#
# After mios-cdi-detect generates /run/cdi/*.{yaml,json} specs for
# every detected GPU vendor, reconcile the AddDevice= drop-ins on
# the background-AI Quadlets (ollama, mios-open-webui, ...) so they
# pick up AMD/Intel iGPUs on the next start. Idempotent: no-op when
# nothing changed; deletes stale drop-ins when an iGPU is removed.
#
# Operator directive 2026-05-17: "ALL background agents work on
# iGPUs detected in the environment!! ... Should be OCI GPU-PV for
# ALL Background AI Agents from the host" + "make sure ALL is in
# code and can deploy Day-0 from clone on any machine".
#
# Lives under /usr/lib/systemd/ (image-immutable) so a fresh Day-0
# install from clone gets the wiring without operator-side action.

[Service]
ExecStartPost=/usr/libexec/mios/mios-gpu-passthrough
